# 10-apple-keyboard.rules
# Reconfigure the keyboard layout of Apple keyboards to match the rest of the
# world's keyboards

SUBSYSTEMS=="input", ATTRS{name}=="Apple Inc. Apple Keyboard", ACTION=="add", RUN+="/usr/bin/fix-apple-keyboard '%E{DEVNAME}'"
