# 10-k810-conf.rules
ACTION=="add", SUBSYSTEM=="hidraw", ATTRS{address}=="FIXME", RUN+="/usr/bin/k810-conf -d $devnode -f on"
