# 20-bt-auto-enable-a2dp.rules
SUBSYSTEM=="bluetooth", ACTION=="add", RUN+="/usr/bin/bt-auto-enable-a2dp"