# 40-a4tech-usb-devices.rules
# A4 Tech Co., Ltd
# 09da:054f
SUBSYSTEM=="usb", ATTRS{idVendor}=="09da", ATTR{idProduct}=="054f", GROUP="users", RUN+="/usr/bin/a4tool init"
