# 50-apple-superdrive.rules
ACTION=="add", ATTRS{product}=="MacBook Air SuperDrive", ATTRS{manufacturer}=="Apple Inc.", ATTRS{removable}=="removable|unknown", DRIVERS=="usb", RUN+="/usr/bin/superdrive-enabler /dev/$kernel"
