# 50-auto-xrandr.rules
SUBSYSTEM=="drm", ACTION=="change", ENV{DISPLAY}=":0", RUN+="/usr/bin/auto-xrandr"
