# 50-spi_ch341_usb.rules
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="1a86", ATTR{idProduct}=="5512", RUN+="/usr/local/bin/spi_ch341_usb_setup.sh"
