# 60-gopro.rules
SUBSYSTEM=="net", ATTRS{idVendor}=="2672", ATTRS{product}=="GoPro HERO9", ACTION=="add", RUN+="/usr/sbin/service gopro_webcam start"
SUBSYSTEM=="net", ATTRS{idVendor}=="2672", ATTRS{product}=="GoPro HERO9", ACTION=="remove", RUN+="/usr/sbin/service gopro_webcam stop"

