# 70-k380.rules
ACTION== "add", KERNEL=="hidraw[0-9]*", RUN+="/usr/bin/k380-function-keys-conf-enable.sh /dev/%k"
