# 80-k380.rules
ACTION== "add", KERNEL=="hidraw[0-9]*", RUN+="/usr/bin/fn_on.sh /dev/%k"
