# 90-acpi_call.rules
SUBSYSTEM=="module", DEVPATH=="/module/acpi_call", RUN+="/bin/sh -c 'find /proc/acpi -name call | xargs /bin/chmod 666'" 
