# 90-lg4l.rules
# This file binds Logitech G keyboard devices to the lg4l modules
# See udev(7) for syntax.
#

# Logitech G19 Gaming Keyboard
SUBSYSTEM=="hid", DRIVER=="generic-usb", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c229", ACTION=="add", RUN+="/bin/sh -c '/sbin/modprobe hid-g19; if test -d /sys/bus/hid/drivers/hid-g19; then /bin/echo -n %k >/sys/bus/hid/drivers/generic-usb/unbind && /bin/echo -n %k >/sys/bus/hid/drivers/hid-g19/bind; fi'"
SUBSYSTEM=="hid", DRIVER=="hid-generic", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c229", ACTION=="add", RUN+="/bin/sh -c '/sbin/modprobe hid-g19; if test -d /sys/bus/hid/drivers/hid-g19; then /bin/echo -n %k >/sys/bus/hid/drivers/hid-generic/unbind && /bin/echo -n %k >/sys/bus/hid/drivers/hid-g19/bind; fi'"

# Logitech G15 Gaming Keyboard (version 1 - Blue)
SUBSYSTEM=="hid", DRIVER=="generic-usb", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c222", ACTION=="add", RUN+="/bin/sh -c '/sbin/modprobe hid-g15; if test -d /sys/bus/hid/drivers/hid-g15; then /bin/echo -n %k >/sys/bus/hid/drivers/generic-usb/unbind && /bin/echo -n %k >/sys/bus/hid/drivers/hid-g15/bind; fi'"
SUBSYSTEM=="hid", DRIVER=="hid-generic", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c222", ACTION=="add", RUN+="/bin/sh -c '/sbin/modprobe hid-g15; if test -d /sys/bus/hid/drivers/hid-g15; then /bin/echo -n %k >/sys/bus/hid/drivers/hid-generic/unbind && /bin/echo -n %k >/sys/bus/hid/drivers/hid-g15/bind; fi'"

# Logitech G15 Gaming Keyboard (version 2 - Orange)
SUBSYSTEM=="hid", DRIVER=="generic-usb", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c227", ACTION=="add", RUN+="/bin/sh -c '/sbin/modprobe hid-g15v2; if test -d /sys/bus/hid/drivers/hid-g15v2; then /bin/echo -n %k >/sys/bus/hid/drivers/generic-usb/unbind && /bin/echo -n %k >/sys/bus/hid/drivers/hid-g15v2/bind; fi'"
SUBSYSTEM=="hid", DRIVER=="hid-generic", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c227", ACTION=="add", RUN+="/bin/sh -c '/sbin/modprobe hid-g15v2; if test -d /sys/bus/hid/drivers/hid-g15v2; then /bin/echo -n %k >/sys/bus/hid/drivers/hid-generic/unbind && /bin/echo -n %k >/sys/bus/hid/drivers/hid-g15v2/bind; fi'"

# Logitech G13 Advanced Gameboard
SUBSYSTEM=="hid", DRIVER=="generic-usb", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c21c", ACTION=="add", RUN+="/bin/sh -c '/sbin/modprobe hid-g13; if test -d /sys/bus/hid/drivers/hid-g13; then /bin/echo -n %k >/sys/bus/hid/drivers/generic-usb/unbind && /bin/echo -n %k >/sys/bus/hid/drivers/hid-g13/bind; fi'"
SUBSYSTEM=="hid", DRIVER=="hid-generic", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c21c", ACTION=="add", RUN+="/bin/sh -c '/sbin/modprobe hid-g13; if test -d /sys/bus/hid/drivers/hid-g13; then /bin/echo -n %k >/sys/bus/hid/drivers/hid-generic/unbind && /bin/echo -n %k >/sys/bus/hid/drivers/hid-g13/bind; fi'"

# Logitech G510 Keyboard
SUBSYSTEM=="hid", DRIVER=="generic-usb", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c22d", ACTION=="add", RUN+="/bin/sh -c '/sbin/modprobe hid-g510; if test -d /sys/bus/hid/drivers/hid-g510; then /bin/echo -n %k >/sys/bus/hid/drivers/generic-usb/unbind && /bin/echo -n %k >/sys/bus/hid/drivers/hid-g510/bind; fi'"
SUBSYSTEM=="hid", DRIVER=="hid-generic", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c22d", ACTION=="add", RUN+="/bin/sh -c '/sbin/modprobe hid-g510; if test -d /sys/bus/hid/drivers/hid-g510; then /bin/echo -n %k >/sys/bus/hid/drivers/hid-generic/unbind && /bin/echo -n %k >/sys/bus/hid/drivers/hid-g510/bind; fi'"
SUBSYSTEM=="hid", DRIVER=="generic-usb", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c22e", ACTION=="add", RUN+="/bin/sh -c '/sbin/modprobe hid-g510; if test -d /sys/bus/hid/drivers/hid-g510; then /bin/echo -n %k >/sys/bus/hid/drivers/generic-usb/unbind && /bin/echo -n %k >/sys/bus/hid/drivers/hid-g510/bind; fi'"
SUBSYSTEM=="hid", DRIVER=="hid-generic", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c22e", ACTION=="add", RUN+="/bin/sh -c '/sbin/modprobe hid-g510; if test -d /sys/bus/hid/drivers/hid-g510; then /bin/echo -n %k >/sys/bus/hid/drivers/hid-generic/unbind && /bin/echo -n %k >/sys/bus/hid/drivers/hid-g510/bind; fi'"

# Logitech G110 Keyboard (no LCD)
SUBSYSTEM=="hid", DRIVER=="generic-usb", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c22b", ACTION=="add", RUN+="/bin/sh -c '/sbin/modprobe hid-g110; if test -d /sys/bus/hid/drivers/hid-g110; then /bin/echo -n %k >/sys/bus/hid/drivers/generic-usb/unbind && /bin/echo -n %k >/sys/bus/hid/drivers/hid-g110/bind; fi'"
SUBSYSTEM=="hid", DRIVER=="hid-generic", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c22b", ACTION=="add", RUN+="/bin/sh -c '/sbin/modprobe hid-g110; if test -d /sys/bus/hid/drivers/hid-g110; then /bin/echo -n %k >/sys/bus/hid/drivers/hid-generic/unbind && /bin/echo -n %k >/sys/bus/hid/drivers/hid-g110/bind; fi'"
