# 99-8bitdo-ultimate-controllers.rules
# dinput
KERNEL=="event[0-9]*", SUBSYSTEM=="input", SUBSYSTEMS=="usb", ENV{ID_VENDOR_ID}=="2dc8", ENV{ID_MODEL_ID}=="3013", TAG+="systemd", ENV{SYSTEMD_WANTS}="8bitdo-ultimate-dinput@$env{DEVNAME}"

# xinput
SUBSYSTEM=="usb", ATTR{idVendor}=="2dc8", ATTR{idProduct}=="3106", ATTR{manufacturer}=="8BitDo", RUN+="/usr/bin/systemctl start 8bitdo-ultimate-xinput@2dc8:3106"

# This device is connected when the above device disconnects
SUBSYSTEM=="usb", ATTR{idVendor}=="2dc8", ATTR{idProduct}=="3109", ATTR{manufacturer}=="8BitDo", RUN+="/usr/bin/systemctl stop 8bitdo-ultimate-xinput@2dc8:3106"
