# 99-k780-config.rules
ATTR{idVendor}=="046d", ATTR{idProduct}=="c52b", RUN+="/usr/bin/k780swap", GROUP="users", MODE="0660"
