# 99-logitech-g923-xbox.rules
# Logitech G923 wheel switch to PC mode
ATTR{idVendor}=="046d", ATTR{idProduct}=="c26d", RUN+="usb_modeswitch '%b/%k'"