# acestream-engine.service
[Unit]
Description=AceStream Engine Service
After=network.target

[Service]
User=acestream
ExecStart=env PYTHON_EGG_CACHE=/var/lib/ACEStream/Python-Eggs /usr/bin/acestream-engine --client-console @/usr/lib/acestream-engine/acestream.conf --live-cache-type memory --state-dir /var/lib/ACEStream --log-file /var/lib/ACEStream/acestream-engine.log

[Install]
WantedBy=multi-user.target 
