# acestream-engine.service
[Unit]
Description=AceStream Engine Service
After=network.target

[Service]
ExecStart=/usr/bin/acestreamengine --client-console --live-cache-type memory --state-dir /usr/share/ACEStream --log-file /usr/share/ACEStream/acestreamengine.log

[Install]
WantedBy=multi-user.target 
