# acme-client@.timer
[Unit]
Description=Check and renew ACME TLS certificates

[Timer]
OnCalendar=daily
RandomizedDelaySec=2h

[Install]
WantedBy=timers.target
