# acme@.timer
[Unit]
Description=Check and renew ACME TLS certificates

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
