# act_runner.service
[Unit]
Description=Act Runner
After=network.target
After=docker.service

[Service]
ExecStart=/usr/bin/act_runner daemon --config /etc/act_runner/act_runner.yaml
ExecReload=/bin/kill -s HUP $MAINPID
WorkingDirectory=/var/lib/act_runner
Type=exec
Restart=always
RestartSec=10
User=act_runner
Environment="HOME=/var/lib/act_runner"

[Install]
WantedBy=multi-user.target