# actkbd.service
[Unit]
Description=Actkbd: Daemon for X-independent shortcuts

[Service]
ExecStart=/usr/bin/actkbd

[Install]
WantedBy=multi-user.target

