# addrwatch.service
[Unit]
Description=Watch Ethernet/IP paring on interface %i
After=systemd-network-wait-online.service
Wants=systemd-network-wait-online.service

[Service]
Type=simple
ExecStart=/usr/bin/addrwatch -q -o /var/lib/addrwatch/%i.dat %i

[Install]
WantedBy=multi-user.target
