# adguardhome.service
[Unit]
Description=AdGuard Home: Network-level blocker
ConditionFileIsExecutable=/var/lib/adguardhome/AdGuardHome
After=syslog.target network-online.target

[Service]
StartLimitInterval=5
StartLimitBurst=10
WorkingDirectory=/var/lib/adguardhome
EnvironmentFile=-/etc/default/adguardhome
ExecStart=/var/lib/adguardhome/AdGuardHome "-s" "run"
StandardOutput=file:/var/log/AdGuardHome.out
StandardError=file:/var/log/AdGuardHome.err
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target