# adsklicensing.service
[Unit]
Description=Autodesk Licensing Service
After=network.target

[Service]
Type=simple
ExecStart=/opt/Autodesk/AdskLicensing/9.2.1.2399/AdskLicensingService/AdskLicensingService --run
Restart=always
RestartSec=0

[Install]
WantedBy=multi-user.target
