# aegir.target
[Unit]
Description=Aegir Hosting System
Wants=aegir.service mailcatcher.service mailcatcher-smtp.socket network.target
Requires=mysqld.service nginx.service php-fpm.service
After=aegir.service mailcatcher.service mailcatcher-smtp.socket mysqld.service network.target nginx.service php-fpm.service

[Install]
WantedBy=multi-user.target
