# aegir.target
[Unit]
Description=Aegir Hosting System
Wants=mailcatcher.service mailcatcher-smtp.socket
After=mailcatcher.service mailcatcher-smtp.socket

[Install]
WantedBy=multi-user.target
