# aegir.service
[Unit]
Description=Aegir queue daemon
Wants=mysqld.service php56-fpm.service
After=mysqld.service php56-fpm.service

[Service]
Type=simple
ExecStart=/usr/bin/drush --quiet @hostmaster hosting-queued
User=http
Restart=always
SuccessExitStatus=1

[Install]
RequiredBy=aegir.target
