# aeond.service
[Unit]
Description=Aeon Full Node
After=network.target

[Service]
User=aeon
Group=aeon
WorkingDirectory=~
RuntimeDirectory=aeon
Type=forking
PIDFile=/run/aeon/aeond.pid
ExecStart=/usr/bin/aeond --config-file /etc/aeond.conf \
    --detach --pidfile /run/aeon/aeond.pid
Restart=always

[Install]
WantedBy=multi-user.target
