# aidecheck.timer
[Unit]
Description=Daily AIDE check

[Timer]
OnCalendar=*-*-* 01:50:00
RandomizedDelaySec=2h
Persistent=true

[Install]
WantedBy=multi-user.target
