# airdcpp.service
[Unit]
Description=AirDC++ Web Client Daemon

[Service]
ExecStart=/usr/bin/airdcppd

[Install]
WantedBy=multi-user.target
