# airdcpp.service
[Unit]
Description=AirDC++ Web Client Daemon
After=network.target

[Service]
Type=simple
User=airdcpp
ExecStart=/usr/bin/airdcppd -c=/var/lib/airdcpp

[Install]
WantedBy=multi-user.target
