# airtime-liquidsoap.service
[Unit]
Description=airtime-liquidsoap daemon
After=local-fs.target remote-fs.target network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/airtime-liquidsoap
ExecReload=/usr/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
