# airtime-media-monitor.service
[Unit]
Description=airtime-media-monitor daemon
After=local-fs.target remote-fs.target network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/airtime-media-monitor
ExecReload=/usr/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
