# aksusbd.service
[Unit]
Description=AKSUSB Daemon
ConditionPathExists=/etc/init.d/aksusbd
After=syslog.target

[Service]
Type=forking
StandardOutput=tty
TimeoutStopSec=10s

ExecStart=/usr/bin/aksusbd
ExecStart=/usr/bin/winehasp
ExecStart=/usr/bin/hasplmd -s

[Install]
WantedBy=multi-user.target
