# aksusbd.service
[Unit]
Description=AKSUSB daemon
Wants=winehasp.service hasplmd.service

[Service]
Type=forking
StandardOutput=tty
TimeoutStopSec=10s

ProtectSystem=full
ProtectHome=true

ExecStart=/usr/bin/aksusbd

[Install]
WantedBy=multi-user.target
Also=winehasp.service
Also=hasplmd.service
