# akula.service
[Unit]
Description=Akula
After=network.target

[Service]
ExecStart=/usr/bin/akula
Restart=on-failure

[Install]
WantedBy=default.target
