# alist.service
[Unit]
Description=alist
After=network.target

[Service]
Type=simple
WorkingDirectory=/usr/bin
ExecStart=/usr/bin/alist server --conf /etc/alist/config.json
Restart=on-failure

[Install]
WantedBy=multi-user.target
