# alist.service
[Unit]
Description=alist
After=network.target

[Service]
Type=simple
User=alist
WorkingDirectory=/var/lib/alist
ExecStart=/usr/bin/alist server --data /var/lib/alist
Restart=on-failure

[Install]
WantedBy=multi-user.target
