# alps-user.service
[Unit]
Description=A simple and extensible webmail
Documentation=https://git.sr.ht/~migadu/alps/tree/master/item/docs
After=network.target
Wants=network-online.target

[Service]
EnvironmentFile=%E/alps/alps.conf
ExecStart=/usr/bin/alps $ALPS_OPTIONS $ALPS_SERVERS
Restart=on-failure

[Install]
WantedBy=default.target
