# amuled@.service
[Unit]
Description=aMule Daemon
After=network.target

[Service]
ExecStart=/usr/bin/amuled

[Install]
WantedBy=multi-user.target
