# ananicy-cpp.openrc
#!/sbin/openrc-run
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

# Init script for ananicy-cpp

pidfile="/var/run/ananicy-cpp.pid"
command="/usr/bin/ananicy-cpp"
start_stop_daemon_args="--wait 30"
command_args=""

stop() {
	ebegin "Stopping ananicy-cpp"
	local KILLOPT="-k"
	while [ -f "${pidfile}" ] && kill -0 $(cat ${pidfile}) ; do
		${command} ${KILLOPT}
		sleep 1
	done
	eend $?
}
