# anbox-container-manager.service
[Unit]
Description=Anbox Container Manager

[Service]
ExecStart=/usr/bin/anbox container-manager --privileged

[Install]
WantedBy=multi-user.target
