# anki-sync-server.service
[Unit]
Description=A personal sync server for anki
After=network.target
StartLimitIntervalSec=30
StartLimitBurst=6
Upholds=nginx.service
PropagatesStopTo=nginx.service

[Service]
Type=simple
User=changeme
Group=changeme
WorkingDirectory=changeme
ExecStart=python -m ankisyncd
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
