# anki-sync-server.service
[Unit]
Description=A sync server for anki
After=network.target
StartLimitIntervalSec=1

[Service]
Type=simple
User=changeme
Group=changeme
WorkingDirectory=changeme
ExecStart=python -m ankisyncd

[Install]
WantedBy=multi-user.target
