# apache-karaf.service
[Unit]
Name=apache-karaf
Description=Apache Karaf application container
After=syslog.target network.target

[Service]
Environment=JAVA_HOME='/usr/lib/jvm/default'
ExecStart=/opt/apache-karaf/bin/karaf server
ExecStop=/opt/apache-karaf/bin/karaf stop

SuccessExitStatus=0 143
Restart=always

LimitNOFILE=102642

SyslogIdentifier=karaf

[Install]
WantedBy=multi-user.target
