# appset-helper.service
[Unit]
Description=Appset Helper Daemon

[Service]
ExecStart=/usr/bin/appset-helper

[Install]
WantedBy=multi-user.target
