# arduino-create-agent.service
[Unit]
Description=arduino create agent

[Service]
Type=Simple
ExecStart=/usr/bin/arduino-create-agent-cli

[Install]
WantedBy=default.target
