# aria2-tracker-update-system.service
[Unit]
Description=Update Aria2c Trackers
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
DynamicUser=yes
User=aria2cd
ExecStart=/usr/bin/update-aria2-tracker --systemd
# Sandboxing
ProtectSystem=strict
ReadOnlyPaths=/
ReadWritePaths=/etc/aria2cd
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectHome=yes
ProtectControlGroups=yes
NoNewPrivileges=yes
