# aria2-tracker-update-user.timer
[Unit]
Description=Update User Aria2c Trackers Daily
PartOf=aria2.target

[Timer]
OnUnitActiveSec=1d

[Install]
WantedBy=aria2.target 
WantedBy=default.target
