# aria2-tracker-update-system.timer
[Unit]
Description=Update System Aria2c Trackers Daily
Before=aria2.service

[Timer]
OnCalendar=8:00
Persistent=yes

[Install]
WantedBy=aria2.target
WantedBy=multi-user.target
