# aria2-system.target
[Unit]
Description=System Aria2 Daemons with Trackers Updated Automatically
After=network-online.target
Wants=network-online.target

