# aria2-tracker-update-user.service
[Unit]
Description=Update Aria2c Trackers
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/update-aria2-tracker
