# aria2-tracker-update-system.service
[Unit]
Description=Update Aria2c Trackers
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/update-aria2-tracker --systemd
RestartSec=60s
Restart=on-failure
# Sandboxing
ProtectSystem=strict
ReadWritePaths=/etc/aria2cd
NoExecPaths=/
ExecPaths=/usr/bin /usr/lib
ProtectHome=yes
NoNewPrivileges=yes
