# aria2-tracker-update-user.timer
[Unit]
Description=Update User Aria2c Trackers Daily
Before=aria2.service

[Timer]
OnUnitActiveSec=1d

[Install]
WantedBy=aria2.target 
WantedBy=default.target
