# aria2-tracker-update-system.timer
[Unit]
Description=Update System Aria2c Trackers Daily
PartOf=aria2.target

[Timer]
OnUnitActiveSec=1d

[Install]
WantedBy=aria2.target
