# aria2-user.target
[Unit]
Description=User Aria2 Daemons with Trackers Updated Automatically
After=network-online.target
Wants=network-online.target 
