# aria2-system.target
[Unit]
Description=System Aria2 Daemons with Trackers Updated Automatically
After=network-online.target
Wants=network-online.target
Requires=aria2.service
After=aria2.service
Wants=aria2-tracker-update.timer
After=aria2-tracker-update.timer

[Install]
WantedBy=multi-user.target
