# aria2-user.service
[Unit]
Description=User Wide Aria2c Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Arguments here will override config file.
ExecStart=/usr/bin/aria2c --daemon=false --enable-rpc=true

[Install]
WantedBy=aria2.target 
WantedBy=default.target
