# asf-user.service
[Unit]
Description=ArchiSteamFarm User Service
Documentation=https://github.com/JustArchi/ArchiSteamFarm/wiki
Wants=network.target

[Service]
ExecStart=/usr/bin/asf --server --path %h/.config/asf

[Install]
WantedBy=default.target
