# asus-charging-threshold.path
[Path]
PathModified=/etc/asus-charging-threshold

[Install]
WantedBy=multi-user.target