# auth-thu.service
[Unit]
Description=Authenticating utility for auth.tsinghua.edu.cn
After=network-online.target

[Service]
ExecStart=auth-thu -D auth -k
ExecStopPost=auth-thu -D deauth
Restart=always

[Install]
WantedBy=default.target
