# auth-thu.timer
[Unit]
Description=Authenticate for auth.tsinghua.edu.cn every 3 minutes

[Timer]
OnBootSec=10
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
