# authdaemond.service
[Unit]
Description=Courier Authentification Daemon

[Service]
Type=forking
User=courier
Group=courier
ExecStart=/usr/bin/authdaemond start
ExecStop=/usr/bin/authdaemond stop
PIDFile=/run/authdaemon/pid
AmbientCapabilities=CAP_DAC_OVERRIDE

[Install]
WantedBy=multi-user.target

