# auth-thu.service
[Unit]
Description=Authenticating utility for auth.tsinghua.edu.cn

[Service]
Type=oneshot
ExecStart=auth-thu --daemonize auth

[Install]
WantedBy=default.target
