# automatic.service
[Unit]
Description=RSS downloader that adds torrents directly to Transmission

[Service]
Type=simple
User=transmission
ExecStart=/usr/bin/automatic -f
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
