# autoupgrade.timer
[Unit]
Description=Automatic Upgrade

[Timer]
OnBootSec=5min
OnActiveSec=20min
OnUnitActiveSec=12h
Unit=autoupgrade.service

[Install]
WantedBy=multi-user.target
