# autopass.service
[Unit]
Description=Socket activated autopass

[Service]
ExecStart=/usr/bin/autopass
