# avantfax.cron.timer
[Unit]
Description=AvantFAX file daily maintenance

[Timer]
OnCalendar=daily
AccuracySec=2h
Persistent=true

[Install]
WantedBy=multi-user.target


