# babeld.service
[Unit]
Description=Babel routing daemon
Documentation=man:babeld(8) http://www.pps.univ-paris-diderot.fr/~jch/software/babel/
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/babeld -D
PIDFile=/var/run/babeld.pid

[Install]
WantedBy=multi-user.target
