# babeld.service
[Unit]
Description=Babel routing daemon
Documentation=man:babeld(8) http://www.pps.univ-paris-diderot.fr/~jch/software/babel/
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/babeld

[Install]
WantedBy=multi-user.target
