# backupRotation.timer
[Unit]
Description=Run backupRotation.service every 24 hours

[Timer]
OnCalendar=*:24/00
