# backup_pacmandb.cron
#!/bin/bash

PACMAN_LOCAL="/var/lib/pacman/local"
BACKUP_DIR="/var/backups/pacman_local_database"

# backup local database
[[ ! -d ${BACKUP_DIR} ]] && mkdir -p ${BACKUP_DIR}
cp -a ${PACMAN_LOCAL} ${BACKUP_DIR}

# create timestamp
cd ${BACKUP_DIR}; echo "$(date)" >> backup.log
