# bacula-fd.service
# systemd service file for Bacula File Deaemon
#
# we do not take the service file provided by bacula package as it needs too
# much customization

[Unit]
Description=Bacula File Deaemon
Requires=network.target
After=network.target remote-fs.target syslog.target

[Service]
Type=forking
Restart=always
ExecStart=/usr/bin/bacula-fd -c /etc/bacula/bacula-fd.conf
IOSchedulingClass=idle
PIDFile=/run/bacula-fd.9102.pid

[Install]
WantedBy=multi-user.target
