# bandwidthd-rotate.service
[Unit]
Description=Daemon for graphing traffic of subnet machines - Rotate

[Service]
ExecStart=/usr/bin/systemctl kill --signal=SIGHUP bandwidthd.service
