# battery-status-collect.timer
[Unit]
Description=Run battery-status-collect every 10 minutes

[Timer]
OnBootSec=10min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
