# battery-notify.service
[Unit]
Description=Notify user on low battery

[Service]
ExecStart=/usr/bin/battery-notify
ProtectSystem=strict
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=default.target
