# batterylife.service
[Unit]
Description=Help user keep battery between 40 and 80%

[Service]
Type=simple
ExecStart=/usr/bin/batterylife
