# battery-status-collect.service
[Unit]
Description=Collect battery statistics

[Service]
Type=oneshot
ExecStart=/usr/bin/battery-status-collect
