# bazarr.service
[Unit]
Description=Bazarr Service
After=network.target

[Service]
User=bazarr
Group=bazarr
Type=simple
ExecStart=/usr/bin/python2 /usr/lib/bazarr/bazarr.py --no-update --config /var/lib/bazarr
Restart=always
TimeoutStopSec=30
SyslogIdentifier=bazarr

[Install]
WantedBy=multi-user.target
