# bbswitch-monitor.service
[Unit]
Description=System tray icon to monitor external GPU on/off status with Bumblebee

[Service]
Type=simple

ExecStart=bbswitch-monitor
Restart=always

[Install]
WantedBy=graphical.target