# beam.service
[Unit]
Description=Beam
After=network.target

[Service]
ExecStart=/usr/bin/beam -headless

[Install]
WantedBy=default.target
