# bedrock.service
[Unit]
Description=Expensify Bedrock Server
After=network.target remote-fs.target nss-lookup.target

[Service]
Type=simple
SyslogIdentifier=bedrock
EnvironmentFile=/etc/bedrock.conf
ExecStart=/usr/bin/bedrock -nodeName $NODE_NAME -db $DATABASE -serverHost $SERVER_HOST -nodeHost $NODE_HOST -priority $PRIORITY -pidfile $PIDFILE -workerThreads $WORKER_THREADS -plugins $PLUGINS -cacheSize $CACHE_SIZE $BEDROCK_EXTRA_ARGS
TimeoutSec=10

[Install]
WantedBy=multi-user.target
