# beehive.service
[Unit]
Description=Beehive

[Service]
User=beehive
ExecStart=/usr/share/webapps/beehive/beehive -config /etc/beehive/beehive.conf

[Install]
WantedBy=multi-user.target
