# bilibilitool-dotnet.service
[Unit]
Description=Bilibili daily task helper
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bilibilitool-dotnet -closeConsoleWhenEnd 1
StandardOutput=null
StandardError=journal


[Install]
WantedBy=multi-user.target
