# bin.service
[Unit]
Description=A paste bin
After=network.target

[Service]
Type=simple
WorkingDirectory=/etc/bin
ExecStart=/usr/bin/bin
Restart=always

[Install]
WantedBy=multi-user.target
