# blackfire-agent.service
[Unit]
Description=Blackfire Agent
After=network.target

[Service]
ExecStart=/usr/bin/blackfire-agent
User=blackfire-agent
Group=blackfire-agent
RuntimeDirectory=blackfire

[Install]
WantedBy=multi-user.target
