# blackd.service
[Unit]
Description=Black Python Formatter Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/blackd

[Install]
WantedBy=default.target
