# blank-html-server.service
[Unit]
Description=Serve blank html pages using node.js

[Service]
ExecStart=/usr/bin/blank-html-server

[Install]
WantedBy=multi-user.target
