# bluealsa.service
[Unit]
Description=BluezAlsa proxy
Requires=bluetooth.service
After=bluetooth.service

[Service]
Type=simple
User=bluealsa
Group=audio
ExecStart=/usr/bin/bluealsa

[Install]
WantedBy=multi-user.target
