# bluealsa.service
[Unit]
Description=BluezAlsa proxy
Requires=bluetooth.service
After=bluetooth.service

[Service]
Type=simple
EnvironmentFile=/etc/conf.d/bluealsa
#User=bluealsa
#Group=audio
ExecStart=/usr/bin/bluealsa $BLUEALSA_CONF

[Install]
WantedBy=multi-user.target
