# bmx7.service
[Unit]
Description=mesh routing network daemon
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/bmx6
ExecStop=/usr/bin/kill -15 $MAINPID
PIDFile=/var/run/bmx6/pid
Restart=on-failure

