# bn-init.service
[Unit]
Description=Bitcoin Node: initialize database in the configured directory
Conflicts=bn.service bs-init.service bs.service

[Service]
Type=oneshot
User=obelisk
Group=obelisk
ExecStart=/usr/bin/bn -i /etc/obelisk/bn/bn.cfg
MemoryDenyWriteExecute=true
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectSystem=full
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
