# boot-esp-sync.service
[Unit]
Description=keep ESP in sync with /boot

[Service]
ExecStart=/usr/bin/boot-esp-sync
