# boot-esp-sync.path
[Unit]
Description=keep ESP in sync with /boot
Requires=efi.mount
After=efi.mount

[Path]
PathChanged=/boot
Unit=boot-esp-sync.service

[Install]
WantedBy=multi-user.target
