# boot-digest-check.service
[Unit]
Description=Check boot partition checksums

[Service]
Type=oneshot
ExecStart=/bin/sh /usr/bin/boot-digest-check

[Install]
WantedBy=multi-user.target
