# brcm43438.service
[Unit]
Description=Broadcom BCM43438 bluetooth HCI
ConditionPathIsDirectory=/proc/device-tree/soc/gpio@7e200000/bt_pins
Before=bluetooth.service
After=dev-ttyAMA0.device

[Service]
Type=simple
ExecStart=/usr/bin/btattach -B /dev/ttyAMA0 -P bcm -S 3000000

[Install]
WantedBy=multi-user.target
