# brscan-skey.service
[Unit]
Description=Brother scan-key-tool

[Service]
User=brscan-skey
Type=forking
ExecStart=/opt/brother/scanner/brscan-skey/brscan-skey
ExecStop=/opt/brother/scanner/brscan-skey/brscan-skey --terminate

[Install]
WantedBy=multi-user.target
