# bs-init.service
[Unit]
Description=Bitcoin Server: initialize database in the configured directory
Conflicts=bs.service bn-init.service bn.service

[Service]
Type=oneshot
User=obelisk
Group=obelisk
ExecStart=/usr/bin/bs -i /etc/obelisk/bs/bs.cfg
MemoryDenyWriteExecute=true
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectSystem=full
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
