# bs.service
[Unit]
Description=Bitcoin Server: a Bitcoin full node and query server
Conflicts=bs-init.service bn-init.service bn.service

[Service]
Type=simple
User=obelisk
Group=obelisk
ExecStart=/usr/bin/bs /etc/obelisk/bs/bs.cfg
ExecReload=/usr/bin/kill -s SIGHUP $MAINPID
ExecStop=/usr/bin/kill -s SIGINT $MAINPID
MemoryDenyWriteExecute=true
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectSystem=full
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
