# btc-cirrus-reporter.timer
[Unit]
Description=Timeline of btc-cirrus-reporter

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
