# btcpayserver.service
[Unit]
Description=Btcpayserver
After=network.target
After=nbxplorer.service

[Service]
Type=oneshot
ExecStart=/usr/bin/btcpayserver
RemainAfterExit=yes
Restart=on-failure
#User=

[Install]
WantedBy=multi-user.target
