# btpd.service
[Unit]
Description=BitTorrent client daemon
After=network-up.target

[Service]
Type=simple
User=btpd
ExecStart=/usr/bin/btpd --no-daemon

[Install]
WantedBy=default.target
