# btrfs-sxbackup@.service
[Unit]
Description=btrfs-sxbackup run on %f

[Service]
Nice=19
IOSchedulingClass=idle
ExecStart=/usr/bin/btrfs-sxbackup run %f
