# btrfs-snap@.service
[Service]
Type=oneshot
ExecStart=/usr/bin/sh -c "/usr/bin/btrfs subvolume snapshot -r %I /.@/%i@$(/usr/bin/date -I)"
