# btrfs-snap-daily@.service
[Service]
Type=oneshot
ExecStart=/usr/bin/btrfs subvolume delete /.@/%i-daily-3
ExecStart=/usr/bin/mv /.@/%i-daily-2 /.@/%i-daily-3
ExecStart=/usr/bin/mv /.@/%i-daily-1 /.@/%i-daily-2
ExecStart=/usr/bin/btrfs subvolume snapshot -r %I /.@/%i-daily-1
