# bubbleupnpserver.service
[Unit]
Description=BubbleUPnP Server
Requires=network-online.target
After=network-online.target


[Service]
User=bubbleupnp
Group=bubbleupnp
WorkingDirectory=/var/lib/bubbleupnp
EnvironmentFile=-/etc/conf.d/bubbleupnpserver
ExecStart=/usr/share/bubbleupnpserver/launch.sh
SuccessExitStatus=1 2 SIGKILL

[Install]
WantedBy=multi-user.target
