# bubbleupnpserver.service
[Unit]
Description=BubbleUPnP Server
Wants=network.target
After=network.target

[Service]
User=bubbleupnp
Group=bubbleupnp
WorkingDirectory=/var/lib/bubbleupnp
EnvironmentFile=-/etc/conf.d/bubbleupnpserver
ExecStart=/usr/share/bubbleupnpserver/launch.sh -dataDir /var/lib/bubbleupnp ${RUN_ARGS}
SuccessExitStatus=1 2 SIGKILL

[Install]
WantedBy=multi-user.target
