# bungeecord-backup.service
[Unit]
Description=BungeeCord Server World Backup
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bungeecord backup
User=bungeecord
Group=bungeecord

[Install]
WantedBy=multi-user.target
