# burp.service
[Unit]
Description=Burp is a backup and restore program.

[Service]
User=root
Type=forking
KillMode=process
ExecStart=/usr/bin/burp -c /etc/burp/burp-server.conf
ExecReload=/bin/kill -HUP $MAINPID
#ExecStop=kill $(head -n1 /var/run/burp.server.pid)

[Install]
WantedBy=multi-user.target
