# burstcoin-wallet.service
[Unit]
Description=BurstCoin Wallet

[Service]
Type=simple
WorkingDirectory=/opt/burstcoin-wallet
ExecStart=/opt/burstcoin-wallet/burst.sh

[Install]
WantedBy=default.target
