# bustd.service
[Unit]
Description = Process killer daemon for out-of-memory scenarios 

[Service]
ExecStart = /usr/bin/bustd -n

[Install]
WantedBy = multi-user.target

