# cachefilesd.service
[Unit]
Description=Local network file caching management daemon

[Service]
Type=forking
ExecStartPre=-/usr/bin/modprobe -qab cachefiles
ExecStart=/usr/bin/cachefilesd -f /etc/cachefilesd.conf

[Install]
WantedBy=multi-user.target
