# caddy.service
[Unit]
Description=A general-purpose HTTP/2 web server
After=network.target nss-lookup.target

[Service]
Type=simple
User=http
Group=http
PrivateTmp=true
LimitNOFILE=8192
ExecStart=/usr/bin/caddy -conf="/etc/caddy/caddy.conf" -host="::" -port="80"
Restart=always

[Install]
WantedBy=multi-user.target
