# caps2esc.service
[Unit]
Description=caps2esc
After=systemd-user-sessions.service

[Service]
ExecStart=/usr/bin/nice -n -20 /usr/bin/caps2esc

[Install]
WantedBy=multi-user.target
