# caps2esc.service
[Unit]
Description=caps2esc
Wants=systemd-udev-settle.service
After=systemd-udev-settle.service

[Service]
ExecStart=/usr/bin/nice -n -20 /usr/bin/caps2esc

[Install]
WantedBy=multi-user.target
