# carbon.service
[Unit]
Description=Graphite/Carbon
After=network.target

[Service]
ExecStart=/usr/bin/carbon-cache.py --nodaemon --config=/etc/carbon/carbon.conf start

[Install]
WantedBy=multi-user.target
