# cardano-node.service
[Unit]
Description=Cardano node
Requires=network.target

[Service]
Type=simple
EnvironmentFile=/etc/conf.d/cardano-node
ExecStart=cardano-node run +RTS -N -RTS \
    --database-path $DB_PATH \
    --socket-path $SOCKET_PATH \
    --host-addr $BIND_IPV4 \
    --host-ipv6-addr $BIND_IPV6 \
    --port $BIND_PORT \
    --config $CONFIG_FILE  \
    --topology $TOPOLOGY_FILE
KillSignal=SIGINT
Restart=always
RestartSec=60
User=cardano
TimeoutStopSec=90

[Install]
WantedBy=multi-user.target
