# catalyst-hook.service
[Unit]
Description=Catalyst's fglrx kernel' module builder
Before=umount.target

[Service]
SysVStartPriority=90
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
ExecStop=/usr/bin/catalyst_build_module check
KillMode=none

[Install]
WantedBy=multi-user.target