# ccpd.service
[Unit]
Description=Canon CAPT daemon
Requires=cups.service
After=cups.service

[Service]
Type=forking
ExecStart=/usr/bin/ccpd

[Install]
WantedBy=printer.target
