# cellframe-node.service
[Unit]
Description=Cellframe Node
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=30
StartLimitBurst=5

[Service]
WorkingDirectory=/opt/cellframe-node
ExecStart=/opt/cellframe-node/bin/cellframe-node
ExecStop=/bin/kill -SIGTERM $MAINPID
Restart=on-failure
LogNamespace=cellframe

[Install]
WantedBy=multi-user.target
