# cf-execd.service
[Unit]
Description=CFEngine Execution Daemon
After=network.target
PartOf=cfengine3.service

[Service]
ExecStart=/var/cfengine/bin/cf-execd
Type=forking
PIDFile=/var/cfengine/cf-execd.pid

[Install]
WantedBy=cfengine3.service
WantedBy=multi-user.target
