# cf-monitord.service
[Unit]
Description=CFEngine Monitoring Daemon
PartOf=cfengine3.service
After=network.target

[Service]
ExecStart=/var/cfengine/bin/cf-monitord
Type=forking
PIDFile=/var/cfengine/cf-monitord.pid

[Install]
WantedBy=cfengine3.service
WantedBy=multi-user.target
