# cf-serverd.service
[Unit]
Description=CFEngine Server Daemon
PartOf=cfengine3.service
After=network.target

[Service]
ExecStart=/var/cfengine/bin/cf-serverd
Type=forking
PIDFile=/var/cfengine/cf-serverd.pid

[Install]
WantedBy=cfengine3.service
WantedBy=multi-user.target
