# cfengine3.service
[Unit]
Description=CFEngine 3 daemons
After=network.target
Requires=cf-execd.service cf-monitord.service cf-serverd.service

[Service]
Type=oneshot
ExecStart=/bin/true
ExecReload=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target