# chann.service
[Unit]
Description=CHANN Anonymous Imageboard 

[Service] 
ExecStart=/usr/bin/startchann
Restart=always 

[Install] 
WantedBy=multi-user.target
