# charybdis.service
[Unit]
Description=Charybdis IRCd
Requires=network.target
After=network.target

[Service]
Type=forking
PIDFile=/run/charybdis/ircd.pid
ExecStart=/usr/bin/charybdis
Restart=always
User=charybdis
Group=charybdis

[Install]
WantedBy=multi-user.target
