# cherokee.service
[Unit]
Description=Cherokee web server
After=network.target

[Service]
ExecStart=/usr/bin/cherokee

[Install]
WantedBy=multi-user.target
