# chinadns@.service
[Unit]
Description=ChinaDNS Service
After=network.target

[Service]
Type=simple
User=nobody
ExecStart=/usr/bin/chinadns -s 114.114.114.114,%i -p 5353 -y 0.3 -d -c /etc/chnroute.txt

[Install]
WantedBy=multi-user.target

