# chromium-fullscreen.service
[Unit]
Description=Systemd service and user to run Chromium in fullscreen mode

[Service]
User=chromium
Group=chromium
PAMName=login
Type=simple
ExecStart=/bin/bash /home/chromium/.systemd/chromium-fullscreen.sh
TimeoutStopSec=30
Restart=always

[Install]
WantedBy=multi-user.target
