# chromium-mpp-setup.service
[Unit]
Description=Setup dummy V4L2 devices for chromium MPP decoding/encoding

[Service]
Type=oneshot
EnvironmentFile=/etc/conf.d/chromium-mpp
ExecStart=/usr/bin/chromium-mpp-setup.sh $PROFILE

[Install]
WantedBy=multi-user.target