# cinder-volume.service
[Unit]
Description=Cinder Volume service

[Service]
ExecStart=/usr/bin/cinder-volume --config-dir /etc/cinder/
User=root

[Install]
WantedBy=multi-user.target
