# cinder-scheduler.service
[Unit]
Description=Cinder Scheduler service

[Service]
ExecStart=/usr/bin/cinder-scheduler --config-dir /etc/cinder/
User=cinder

[Install]
WantedBy=multi-user.target
