# clamfs.service
[Unit]
Description=ClamFS

[Service]
Type=forking
ExecStart=/usr/bin/clamfs /etc/clamav/clamfs.xml

[Install]
WantedBy=multi-user.target
