# clashr@.service
[Unit]
Description=A rule based proxy in Go with SSR protocol support for %i.
After=network.target

[Service]
Type=simple
User=%i
Restart=on-abort
ExecStart=/usr/bin/clashr

[Install]
WantedBy=multi-user.target
