# clash-meta.service
[Unit]
Description=A rule based proxy in Go.
After=network.target

[Service]
Type=exec
Restart=on-abort
ExecStart=/usr/bin/clash-meta -d /etc/clash-meta

[Install]
WantedBy=default.target
