# clickhouse-server.service
[Unit]
Description=ClickHouse server

[Service]
Type=simple
PIDFile=/var/run/clickhouse-server.pid
ExecStart=/usr/sbin/clickhouse-server
WorkingDirectory=/etc/clickhouse-server
User=clickhouse
Group=clickhouse

[Install]
WantedBy=multi-user.target
