# cloud9@.service
[Unit]
Description=Cloud9 Ide On Custom Port

[Service]
Type=simple
User=cloud9
ExecStart=/usr/bin/node /opt/cloud9/server.js -w /opt/cloud9/%I -p %I -l 0.0.0.0 -a :

[Install]
WantedBy=multi-user.target
