# code-server.service
[Unit]
Description=code-server user service

[Service]
Type=Simple
ExecStart=/usr/bin/code-server

[Install]
WantedBy=default.target
