# codimd.service
[Unit]
Description=CodiMD real-time collaborative markdown editor
Documentation=file:/opt/codimd/docs/
After=network.target


[Service]
Type=simple

Environment=NODE_ENV=production
User=codimd
Group=codimd
WorkingDirectory=/usr/share/webapps/codimd
ExecStart=/usr/bin/node app.js

StandardOutput=journal
StandardError=inherit
SyslogIdentifier=codimd


[Install]
WantedBy=multi-user.target
