# conduit.service
[Unit]
Description=Conduit Matrix Homeserver
After=network.target

[Service]
EnvironmentFile=-/etc/conduit/config
User=conduit
Group=conduit
Type=simple
Restart=always
ExecStart=/usr/bin/conduit
PrivateTmp=true

[Install]
WantedBy=multi-user.target