# confluence-companion-app-for-linux.service
[Unit]
Description=Confluence Companion App for Linux

[Service]
Environment="DISPLAY=:0"
Environment="XAUTHORITY=%h/.Xauthority"
ExecStart=/usr/bin/companion.py

[Install]
WantedBy=graphical.target
