# confluence-companion.service
[Unit]
Description=Confluence Companion (Legacy Version)

[Service]
Environment="DISPLAY=:0"
Environment="XAUTHORITY=%h/.Xauthority"
ExecStart=/usr/bin/companion.py

[Install]
WantedBy=graphical.target
