# consul-template.target
[Unit]
Description=Target file for consul-template@ service templates

[Install]
WantedBy=multi-user.target
