# controlloid-server.service
[Unit]
Description=The server for the Android app Controlloid

[Service]
ExecStart=/usr/bin/controlloid-server

[Install]
WantedBy=default.target