# continuous-silence.service
[Unit]
Description=continuous-silence

[Service]
ExecStart=/usr/bin/play -qn

[Install]
WantedBy=default.target
