# core-daemon.service
[Unit]
Description=Common Open Research Emulator Service
After=network.target

[Service]
Type=simple
PIDFile=/run/cored.pid
ExecStart=/usr/bin/python2 /usr/sbin/core-daemon -i /run/cored.pid

[Install]
WantedBy=multi-user.target
