# coreos-cloudinit.service
[Unit]
Description=Cloudinit from EC2-style metadata
[Service]
Environment=PATH=${PATH}:/usr/bin
Type=oneshot
ExecStart=/usr/bin/coreos-cloudinit --oem=ec2-compat
