# courier-esmtpd-msa.service
[Unit]
Description=Courier ESMTP-MSA Daemon
Requires=authdaemond.service
After=authdaemond.service

[Service]
Type=forking
ExecStart=/usr/sbin/esmtpd-msa start
ExecStop=/usr/sbin/esmtpd-msa stop
ExecReload=/usr/sbin/esmtpd-msa restart
PIDFile=/run/courier/esmtpd-msa.pid

[Install]
WantedBy=multi-user.target
