# courier-imapd-ssl.service
[Unit]
Description=Courier IMAP Daemon (SSL)
Requires=authdaemond.service
After=authdaemond.service

[Service]
Type=forking
ExecStart=/usr/sbin/imapd-ssl start
ExecStop=/usr/sbin/imapd-ssl stop
ExecReload=/usr/sbin/imapd-ssl restart
PIDFile=/run/courier/imapd-ssl.pid

[Install]
WantedBy=multi-user.target
