# courier-esmtpd-ssl.service
[Unit]
Description=Courier ESMTP Daemon
Requires=authdaemond.service
After=authdaemond.service courier-courierfilter.service

[Service]
Type=forking
ExecStart=/usr/sbin/esmtpd-ssl start
ExecStop=/usr/sbin/esmtpd-ssl stop
ExecReload=/usr/sbin/esmtpd-ssl restart
PIDFile=/run/courier/esmtpd-ssl.pid

[Install]
WantedBy=multi-user.target
