# courier-esmtpd.service
[Unit]
Description=Courier ESMTP Daemon
Requires=authdaemond.service
After=authdaemond.service, courier-courierfilter.service

[Service]
Type=forking
ExecStartPre=[ -d /var/run/courier ] || mkdir -p /var/run/courier && chown courier:courier /var/run/courier
ExecStart=/usr/sbin/esmtpd start
ExecStop=/usr/sbin/esmtpd stop
PIDFile=/var/run/courier/esmtpd.pid

[Install]
WantedBy=multi-user.target
