# courier-pop3d-ssl.service
[Unit]
Description=Courier POP3 (SSL) Daemon
Requires=authdaemond.service
After=authdaemond.service

[Service]
Type=forking
EnvironmentFile=/etc/courier-imap/pop3d
ExecStart=/usr/lib/courier-imap/pop3d-ssl.rc start
ExecStop=/usr/lib/courier-imap/pop3d-ssl.rc stop
PIDFile=/var/run/courier/pop3d-ssl.pid

[Install]
WantedBy=multi-user.target
