# courier-mkdhparams.service
[Unit]
Description=Generate Diffie Hellman parameters for Courier MTA

[Service]
Type=oneshot
ExecStart=/usr/bin/mkdhparams
