# courier-pop3d-ssl.service
[Unit]
Description=Courier POP3 (SSL) Daemon
Requires=authdaemond.service
After=authdaemond.service

[Service]
Type=forking
User=courier
Group=courier
EnvironmentFile=/etc/courier-imap/pop3d
ExecStart=/usr/lib/courier-imap/pop3d-ssl.rc start
ExecStop=/usr/lib/courier-imap/pop3d-ssl.rc stop
PIDFile=/run/courier/pop3d-ssl.pid
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_SETGID CAP_SETUID

[Install]
WantedBy=multi-user.target
