# courier.service
[Unit]
Description=Courier Daemon
Requires=authdaemond.service
After=authdaemond.service

[Service]
Type=forking
ExecStartPre=[ -d /var/run/courier ] || mkdir -p /var/run/courier && chown courier:courier /var/run/courier
ExecStart=/usr/sbin/courier start
ExecStop=/usr/sbin/courier stop

[Install]
WantedBy=multi-user.target
