# courier.service
[Unit]
Description=Courier Daemon
Requires=authdaemond.service
After=authdaemond.service

[Service]
Type=forking
ExecStart=/usr/sbin/courier start
ExecStop=/usr/sbin/courier stop

[Install]
WantedBy=multi-user.target
