# cowrie.service
[Unit]
Description=Cowrie Honey Pot
After=network.target

[Service]
Type=forking
User=cowrie
Group=cowrie
WorkingDirectory=/opt/cowrie
PIDFile=/opt/cowrie/cowrie.pid
ExecStart=/usr/bin/twistd -l log/cowrie.log --pidfile cowrie.pid cowrie

[Install]
WantedBy=multi-user.target
