# cpcache.service
[Unit]
Description=central pacman cache
After=network.target

[Service]
Type=simple
# Shutdown can be somewhat sluggish on low-end devices like the RPi.
TimeoutSec=30
SuccessExitStatus=143
KillMode=process
User=cpcache
Group=cpcache
Environment="ELIXIR_ERL_OPTIONS=-sname cpcache -mnesia dir '/var/cache/cpcache/mnesia'"
ExecStart=/usr/share/cpcache/bin/cpcache foreground

[Install]
WantedBy=multi-user.target
