# cpcache.service
[Unit]
Description=central pacman cache
After=network.target

[Service]
Type=simple
KillMode=process
TimeoutSec=10
User=cpcache
Group=cpcache
ExecStart=/usr/share/cpcache/bin/cpcache foreground

[Install]
WantedBy=multi-user.target
