# cpcache.service
[Unit]
Description=central pacman cache
After=network.target

[Service]
Type=simple
KillMode=process
TimeoutSec=10
User=cpcache
Group=cpcache
Environment="ELIXIR_ERL_OPTIONS=-sname cpcache -mnesia dir '/var/cache/cpcache/mnesia'"
ExecStart=/usr/share/cpcache/bin/cpcache foreground

[Install]
WantedBy=multi-user.target
