# cpufreqd.service
[Unit]
Description=Cpufreqd daemon
After=acpid.service

[Service]
ExecStart=/usr/bin/cpufreqd -D
ExecStop=/bin/kill -INT $MAINPID

[Install]
WantedBy=multi-user.target
