# craftbukkit.service
[Unit]
Description=Minecraft server

[Service]
User=craftbukkit
Group=craftbukkit
RuntimeDirectory=craftbukkit
WorkingDirectory=/srv/craftbukkit
Environment=JAVA_ARGS=-Xmx1G RUNTIME_DIR=/run/craftbukkit
ExecStartPre=/usr/share/craftbukkit-spigot/init-rcon.sh
ExecStart=/usr/bin/java $JAVA_ARGS -jar /usr/share/java/craftbukkit-spigot/%p.jar --noconsole
ExecStop=/usr/share/craftbukkit-spigot/notify-shutdown.sh
PrivateTmp=true
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
