# craftbukkit@.service
[Unit]
Description=Minecraft server for %f

[Service]
User=craftbukkit
Group=craftbukkit
RuntimeDirectory=craftbukkit@%i
WorkingDirectory=%f
Environment=JAVA_ARGS=-Xmx1G RUNTIME_DIR=/run/craftbukkit@%i
ExecStartPre=/usr/share/craftbukkit-spigot/init-rcon.sh
ExecStart=/usr/bin/java $JAVA_ARGS -jar /usr/share/java/craftbukkit-spigot/%p.jar --noconsole
ExecStop=/usr/share/craftbukkit-spigot/notify-shutdown.sh
PrivateTmp=true
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
