# crashplan.service
[Unit]
Description=CrashPlan Backup Engine
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
PIDFile=/opt/crashplan/CrashPlanEngine.pid
EnvironmentFile=/opt/crashplan/bin/run.conf
WorkingDirectory=/opt/crashplan
ExecStart=/opt/crashplan/bin/CrashPlanEngine start
ExecStop=/opt/crashplan/bin/CrashPlanEngine stop
Restart=always

[Install]
WantedBy=multi-user.target
