# crashplan-pro_upgrade.service
[Unit]
Description=CrashPlanPRO Backup Engine Upgrade
After=network.target
Requires=crashplan-pro.service
PropagatesReloadTo=crashplan-pro.service
OnSuccess=crashplan-pro.service
RefuseManualStart=True

[Service]
Type=oneshot
EnvironmentFile=/opt/crashplan/bin/run.conf
WorkingDirectory=/opt/crashplan
ExecStart=/opt/crashplan/bin/upgrade.sh
Restart=no
StandardOutput=append:/opt/crashplan/log/service_upgrade.log
StandardError=append:/opt/crashplan/log/service_upgrade.log
