# crashplan-pro.service
[Unit]
Description=CrashPlanPRO Backup Engine
After=network.target

[Service]

Type=forking
PIDFile=/opt/crashplan/Code42Service.pid
EnvironmentFile=/opt/crashplan/bin/run.conf
WorkingDirectory=/opt/crashplan
ExecStart=/opt/crashplan/bin/service.sh start
ExecStop=/opt/crashplan/bin/service.sh stop
Restart=always

[Install]
WantedBy=multi-user.target
