# crashplan-pro.service
[Unit]
Description=CrashPlanPRO Backup Engine
After=network.target

[Service]

Type=forking
PIDFile=/opt/crashplan-pro/CrashPlanEngine.pid
EnvironmentFile=/opt/crashplan-pro/bin/run.conf

WorkingDirectory=/opt/crashplan-pro

ExecStart=/opt/crashplan-pro/bin/CrashPlanEngine start
ExecStop=/opt/crashplan-pro/bin/CrashPlanEngine stop

[Install]
WantedBy=multi-user.target
