# crashplan-pro_upgrade.path
[Unit]
Description="File monitor for updates from the Crashplan Daemon"
Requires=crashplan-pro.service
After=crashplan-pro.service


[Path]
PathExistsGlob=/opt/crashplan/upgrade/*/install.sh
Unit=crashplan-pro_upgrade.service

[Install]
WantedBy=multi-user.target
