# crate.service
[Unit]
Description=CrateDB Server
After=network.target

[Service]
User=crate
EnvironmentFile=-/etc/conf.d/crate
ExecStart=/usr/share/crate/bin/crate -Des.path.conf=/etc/crate -Des.path.home=/usr/share/crate -Des.path.logs=/var/log/crate -Des.path.work=/var/lib/crate

LimitMEMLOCK=infinity
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
