# crate.service
[Unit]
Description=CrateDB Server
After=network.target

[Service]
User=crate
EnvironmentFile=-/etc/conf.d/crate
ExecStart=/usr/share/crate/bin/crate -Cpath.conf=/etc/crate -Cpath.home=/usr/share/crate -Cpath.logs=/var/log/crate -Cpath.work=/var/lib/crate

LimitMEMLOCK=infinity
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
