# crestic@.timer
[Unit]
Description=Daily crestic %I backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target
