# crestic-forget@.service
[Unit]
Description=crestic %I forget

[Service]
KillSignal=SIGINT
ExecStart=/usr/bin/crestic %I forget
