# crestic@.service
[Unit]
Description=crestic backup %I

[Service]
Environment="CRESTIC_CONFIG_FILE=/etc/crestic/crestic.cfg"
ExecStart=/usr/bin/crestic %I backup
