# crestic@.timer
[Unit]
Description=Daily crestic backup %I

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target
