# crestic@.service
[Unit]
Description=crestic %I backup

[Service]
ExecStart=/usr/bin/crestic %I backup
