# croodle-cron.timer
[Unit]
Description=Croodle daily deletion of outdated polls timer

[Timer]
OnCalendar=daily
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target