# crowdsec-upgrade.service
[Service]
Type=oneshot
ExecStart=/usr/bin/cscli --error hub upgrade
ExecStartPost=systemctl reload crowdsec
