# cryptpad.service
[Unit]
Description=CryptPad service

[Service]
ExecStart=/usr/bin/node /usr/share/cryptpad/server.js
WorkingDirectory=/usr/share/cryptpad
User=cryptpad

[Install]
WantedBy=multi-user.target

