# csync2.socket
[Socket]
ListenStream=30865
Accept=yes

[Install]
WantedBy=sockets.target

