# cuberite.service
[Unit]
Description=Cuberite Server
After=local-fs.target network.target

[Service]
Type=forking
ExecStart=/usr/bin/cuberite start
ExecStop=/usr/bin/cuberite stop
User=cuberite
Group=cuberite

[Install]
WantedBy=multi-user.target
