# cuberite-backup.service
[Unit]
Description=Cuberite Server World, Plugins and Configuration-File Backup
After=local-fs.target

[Service]
Type=forking
ExecStart=/usr/bin/cuberite backup
User=cuberite
Group=cuberite

[Install]
WantedBy=multi-user.target
