# curly-turtle@.service
[Unit]
Description=Curly Turtle
After=network.target

[Service]
User=curly-turtle
ExecStart=/usr/bin/curly-turtle -c /etc/curly-turtle/%i.yml

[Install]
WantedBy=multi-user.target

