# cvsync.service
[Unit]
Description=Synchronize CVS trees
RefuseManualStart=no
RefuseManualStop=yes
 
[Service]
User=cvsync
Type=oneshot
ExecStart=/usr/bin/cvsync
