# dam.service
[Unit]
Description=Dam
Requires=network.target

[Service]
Type=simple
WorkingDirectory=/usr/lib/dam
ExecStart=/bin/bash 'bundle exec rackup -o 0.0.0.0 -p 9669'
TimeoutSec=30
RestartSec=15s
Restart=always

[Install]
WantedBy=multi-user.target
