# darkice@.service
[Unit]
Description=DarkIce Live Stream %i
Documentation=man:darkice(1)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/darkice -c %h/.config/darkice/%i.conf

[Install]
WantedBy=default.target
