# darkstar-game.service
[Unit]
Description=Final Fantasy XI server emulator (game)
After=network.target nss-lookup.target

[Service]
User=darkstar
Group=darkstar
PIDFile=/run/darkstar-game.pid
WorkingDirectory=/var/lib/darkstar
ExecStart=/usr/bin/dsgame
Restart=always

[Install]
WantedBy=multi-user.target
