# dashd@.service
[Unit]
Description=Dash Daemon
After=network.target 

[Service]
User=%i
Group=users
Type=simple
ExecStart=/usr/bin/dashd

[Install]
WantedBy=multi-user.target
