# dcssd.service
[Unit]
Description=Server for injecting CSS into chrome/chromium
StopWhenUnneeded=true

[Service]
ExecStart=/usr/bin/dcssd
