# ddh.service
[Unit]
Description=DynamicDisplayHandler
After=multi-user.target

[Service]
User=riddlerxenon
Environment=DISPLAY=:0
ExecStart=/ -c 
Restart=always

[Install]
WantedBy=multi-user.target
