# ddns-watchdog-server.service
[Unit]
Description=ddns-watchdog-server Service
After=network.target

[Service]
Type=simple
ExecStart=ddns-watchdog-server -c /etc/ddns-watchdog/
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
