# ddns-watchdog-server.service
[Unit]
Description=watchdog-ddns-server Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/watchdog-ddns-server -conf_path /etc/watchdog-ddns/
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
