# ddns-watchdog-client.service
[Unit]
Description=ddns-watchdog-client Service
After=network.target

[Service]
Type=simple
ExecStart=ddns-watchdog-client -c /etc/ddns-watchdog/
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
