# ddns-go.service
[Unit]
Description=简单好用的DDNS。自动更新域名解析到公网IP(支持阿里云、腾讯云dnspod、Cloudflare、华为云)

[Service]
StartLimitInterval=5
StartLimitBurst=10
ExecStart=/usr/bin/ddns-go "-l" ":9876" "-f" "300" "-c" "/etc/ddns-go.yaml"
Restart=always
RestartSec=120

[Install]
WantedBy=multi-user.target
