# ddns-watchdog-client.service
[Unit]
Description=watchdog-ddns-client Service
After=network.target

[Service]
Type=simple
ExecStart=watchdog-ddns-client -c /etc/watchdog-ddns/
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
