# deepl-api.service
[Unit]
Description=Free deepL api, written in rust.
After=network.target

[Service]
Type=simple
ExecStart=deepl-api
Restart=on-abort
RestartSec=10

[Install]
WantedBy=multi-user.target