# dell-g5se-fanctl.service
[Unit]
Description=Control fans for Dell G5SE laptops.

[Service]
Type=simple
StandardOutput=null

Restart=on-failure
RestartSec=10

EnvironmentFile=/etc/dell-g5se-fanctl.conf

ExecStart=dell-g5se-fanctl -l -ut $DELL_TEMPS -t $DELL_INTERVAL $DELL_FLAGS

[Install]
WantedBy=multi-user.target
