# deluged.service
[Unit]
Description=Deluge Bittorrent Client Daemon
Documentation=man:deluged
After=network-online.target

[Service]
User=deluge
EnvironmentFile=/etc/conf.d/deluged
ExecStart=/usr/bin/deluged -P /run/deluge/deluged.pid $OPTIONS
Type=forking
PIDFile=/run/deluge/deluged.pid
UMask=007
Restart=on-failure
TimeoutStopSec=300

[Install]
WantedBy=multi-user.target
