# deluged.service
[Unit]
Description=Deluge Bittorrent Client Daemon
Documentation=man:deluged(1)
After=network-online.target

[Service]
EnvironmentFile=/etc/conf.d/deluged
User=deluge
EnvironmentFile=/etc/conf.d/deluged
ExecStart=/usr/bin/deluged -d $OPTIONS
UMask=007
Restart=on-failure
TimeoutStopSec=300

[Install]
WantedBy=multi-user.target
