# deluge-web.service
[Unit]
Description=Deluge Bittorrent Client WebUI
Documentation=man:deluge-web(1)
After=online.target deluged.service
Wants=deluged.service

[Service]
Umask=027
User=deluge
EnvironmentFile=/etc/conf.d/deluged
ExecStart=/usr/bin/deluge-web -d $WEB_OPTIONS
Restart=on-failure

[Install]
WantedBy=multi-user.target
