# desired@.service
[Unit]
Description=Desire Daemon
After=network.target

[Service]
User=%i
Group=users
Type=simple
ExecStart=/usr/bin/desired

[Install]
WantedBy=multi-user.target
