# dev-binderfs.mount
[Mount]
Where=/dev/binderfs
What=binder
Type=binder
Options=nofail
