# devpi-server.service
[Unit]
Description=Devpi Server
Requires=network-online.target
After=network-online.target

[Service]
Restart=on-success
ExecStart=/usr/bin/devpi-server --serverdir /var/lib/devpi
User=devpi

[Install]
WantedBy=multi-user.service
