# dgraph-zero.service
[Unit]
Description=Dgraph database

[Service]
ExecStart=/usr/bin/dgraph zero --port_offset -2000
Restart=on-failure
RestartSec=1

[Install]
WantedBy=multi-user.target

