# diamond.service
[Unit]
Description=Diamond daemon
After=network.target

[Service]
ExecStart=/usr/bin/diamond --foreground

[Install]
WantedBy=multi-user.target
