# dinit.service
[Unit]
Description=dinit - Service monitoring / "init" system

[Service]
ExecStart=/usr/bin/dinit --system

[Install]
WantedBy=default.target

