# dinit-user.service
[Unit]
Description=dinit - Service monitoring / "init" system

[Service]
ExecStart=/usr/bin/dinit --user

[Install]
WantedBy=default.target

