# discourse-kill.timer
[Unit]
Description=Trigger for daily discourse RAM cleanup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timer.target
