# dispad.service
[Unit]
Description=dispad daemon

[Service]
ExecStart=/usr/bin/dispad -F

[Install]
WantedBy=user-applications.target
