# dn42-peerfinder-client.timer
[Unit]
Description=Run dn42 PeerFinder client

[Timer]
OnBootSec=10min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
