# dnscrypt-proxy.service
[Unit]
Description=DNSCrypt client proxy
Documentation=https://github.com/jedisct1/dnscrypt-proxy
After=network.target
Before=nss-lookup.target

[Service]
Type=simple
ExecStart=/usr/bin/dnscrypt-proxy

[Install]
WantedBy=multi-user.target
