# dnsproxy-adguard@.service
[Unit]
Description=Simple DNS proxy with DoH, DoT, and DNSCrypt support by AdguardTeam
After=network.target

[Service]
Type=simple
User=nobody
Group=nobody
ExecStart=/usr/bin/dnsproxy-adguard -l 127.0.0.1 -p %i -u tls://8.8.8.8 -u tls://8.8.4.4 -u https://mozilla.cloudflare-dns.com/dns-query -u https://dns.adguard.com/dns-query -u tls://9.9.9.9 -u tls://dns.adguard.com --cache --all-servers 

[Install]
WantedBy=default.target
