# dnsproxy-adguard.service
[Unit]
Description=Simple DNS proxy with DoH, DoT, and DNSCrypt support by AdguardTeam
After=network.target

[Service]
Type=simple
User=nobody
Group=nobody
EnvironmentFile=/etc/conf.d/dnsproxy-adguard
ExecStart=/usr/bin/dnsproxy-adguard -l $ADDRESS -p $PORT $UPSTREAMS $OTHER_PARAMS 

[Install]
WantedBy=default.target
