# docker-dnsmasq.service
[Unit]
Description=Daemon to add docker containers to dnsmasq
Requires=docker.service
After=dnsmasq.service

[Service]
Type=idle
ExecStart=/usr/bin/docker-dnsmasq daemon
Restart=always

[Install]
WantedBy=multi-user.target
