# docker-prune.service
[Unit]
Description=Remove unused docker data

[Service]
Type=oneshot
ExecStart=/usr/bin/docker system prune -af --filter "until=240h"
