# docker_makepkg.install
post_install() {
        systemctl start docker-makepkg
}

post_remove() {
        systemctl disable docker-makepkg.timer
        systemctl stop docker-makepkg
        docker image rm makepkg
}
