# docker-prune.timer
[Unit]
Description=Run docker-prune daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
