# dogcom-p.service
[Unit]
Description=dogcom
Requires=network.target
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/dogcom -m pppoe -c /etc/dogcom.d/pppoe.conf -v

[Install]
WantedBy=multi-user.target
