# dom4-server.service
[Unit]
Description=Dominions 4 Headless Server
After=network.target

[Service]
ExecStart=/usr/bin/dom4-server
Type=simple
User=dom4
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
