# domjudge-submitd.service
[Unit]
Description=Domjudge submit daemon
After=network.target

[Service]
ExecStart=/usr/bin/submitdaemon
Restart=on-failure
KillSignal=SIGINT
SendSIGKILL=no

[Install]
WantedBy=multi-user.target

