# domjudge-judgehostd.service
[Unit]
Description=Domjudge judgehost daemon
After=network.target

[Service]
ExecStart=/usr/bin/judgedaemon
Restart=on-failure
User=domjudge
KillSignal=SIGINT
SendSIGKILL=no

[Install]
WantedBy=multi-user.target

