# domoticz.service
[Unit]
Description=Domoticz Daemon
After=network-online.target

[Service]
ExecStart=/opt/domoticz/domoticz

[Install]
WantedBy=multi-user.target
