# domoticz.service
[Unit]
Description=Domoticz Daemon
After=network.target

[Service]
User=http
ExecStart=/opt/domoticz/domoticz -www 8080 -daemon
WorkingDirectory=/opt/domoticz
RestartSec=5
Restart=on-failure

[Install]
WantedBy=multi-user.target
