# dowse.service
[Unit]
Description=Dowse
After=network.target

[Service]
User=dowse
ExecStart=/usr/bin/dowse start

[Install]
WantedBy=multi-user.target
