# drone.service
[Unit]
Description=Continuous Integration platform 
After=network.target

[Service]
User=drone
Type=simple
ExecStart=/usr/bin/droned --config /etc/drone/drone.toml
ExecReload=/usr/bin/kill -HUP $MAINPID
KillSignal=SIGINT
LimitNOFILE=8192
PrivateDevices=yes

[Install]
WantedBy=multi-user.target
