# dropbear-socket@.service
[Unit]
Description=dropbear SSH service (socket-activated)

[Service]
ExecStart=/usr/bin/dropbear -iwgRT 1
StandardInput=socket
