# dropbear-socket@.service
[Unit]
Description=dropbear SSH service (socket-activated)

[Service]
ExecStart=/usr/bin/dropbear -iwgRT1 -K180
StandardInput=socket
