# duckdns.service
[Unit]
Description=Update DuckDNS domain addresses

[Service]
Type=oneshot
ExecStart=/usr/bin/duckdns

[Install]
WantedBy=network-online.target
