# duckdns.timer
[Unit]
Description=Run DuckDNS every 5 minutes

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
