# duckdns.service
[Unit]
Description=DuckDNS update job

[Install]
WantedBy=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/bin/duckdns
