# dump1090.service
[Unit]
Description=RtlSdr mode S decoder
Wants=network-online.target
After=syslog.target network-online.target lighttpd.service

[Service]
Restart=always
RestartSec=30
RuntimeDirectory=dump1090
ExecStart=/usr/bin/dump1090 --quiet --net --gain -10 --write-json /var/run/dump1090

[Install]
WantedBy=multi-user.target
