# dynafire.service
[Unit]
Description=dynafire - real-time threat detection for any Linux system
After=network.target

[Install]
WantedBy=multi-user.target

[Service]
User=root
Type=simple
ExecStart="/usr/bin/dynafire"
TimeoutStopSec=20
KillMode=process
Restart=on-failure
