# dynamic-dns-netcup-api.timer
[Unit]
Description=periodic update of netcup dynDNS
After=network-online.target network.target

[Timer]
OnBootSec=30s
OnUnitActiveSec=5m
#Persistent=true
Unit=dynamic-dns-netcup-api.service

[Install]
WantedBy=timers.target
