# dyndnsc.service
[Unit]
Description=dyndnsc
After=network.target

[Service]
ExecStart=/usr/bin/dyndnsc --loop --config /etc/dyndnsc.conf

[Install]
WantedBy=multi-user.target
