# e2guardian.service
[Unit]
Description=E2guardian web filtering
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/e2guardian

[Install]
WantedBy=multi-user.target
