# easyeda-router.service
[Unit]
Description=EasyEDA Local Auto Router Server
After=network.target

[Service]
Type=simple
DynamicUser=yes
ExecStart=/usr/bin/easyeda-router
Restart=always

[Install]
WantedBy=multi-user.target