# ecbd.service
[Unit]
Description=Epson Printer Utility Daemon

[Service]
Type=forking
ExecStart=/opt/epson-printer-utility/lib/epson-backend/ecbd
Restart=on-failure

[Install]
WantedBy=multi-user.target
