# eGTouchD.service
[Unit]
Description=eGalax touchscreen device driver

[Service]
ExecStart=/usr/bin/eGTouchD
Type=forking

[Install]
WantedBy=graphical.target
