# ehforwarderbot@.service
[Unit]
Description=EH Forwarder Bot Service
Documentation=https://ehforwarderbot.readthedocs.io/zh_CN/latest/
Wants=network-online.target systemd-networkd-wait-online.service
After=network-online.target

[Service]
Type=simple
User=%i
ExecStart=/usr/bin/ehforwarderbot
Restart=on-failure
KillSignal=SIGINT
TimeoutStopSec=10s

[Install]
WantedBy=multi-user.target