# eibd.socket
[Unit]
Conflicts=eibd.service

[Socket]
ListenDatagram=6720

[Install]
WantedBy=sockets.target
