# elasticsearch-keystore.service
[Unit]
Description=Elasticsearch Keystore Generation
ConditionPathExists=|!/etc/elasticsearch6/elasticsearch.keystore

[Service]
Type=oneshot
Group=elasticsearch
UMask=0007
ExecStart=/usr/share/elasticsearch6/bin/elasticsearch-keystore create
RemainAfterExit=yes
