# elasticsearch-keystore@.service
[Unit]
Description=Elasticsearch Keystore Generation
ConditionPathExists=|!/etc/elasticsearch6/%I/elasticsearch.keystore

[Service]
Type=oneshot
Group=elasticsearch
UMask=0007
Environment=ES_PATH_CONF=/etc/elasticsearch6/%I
ExecStart=/usr/share/elasticsearch6/bin/elasticsearch-keystore create
RemainAfterExit=yes
