# endpoint-verification.service
[Unit]
Description=Endpoint Verification Service

[Service]
Type=oneshot
ExecStart=/opt/google/endpoint_verification/endpoint-verification.sh

[Install]
WantedBy=graphical.target
