# entrance.service
[Unit]
Description=Enlightenment Display Manager
After=systemd-user-sessions.service

[Service]
ExecStart=/usr/bin/entrance --nodaemon
Restart=always

[Install]
Alias=display-manager.service
