# epgd.service
[Unit]
Description = Database driven EPG Data collector
After = mariadb.service

[Service]
Type = forking
ExecStart = /usr/bin/epgd
Restart = always

[Install]
WantedBy = multi-user.target
