# epoptes-client.service
[Unit]
Description=Epoptes-client system daemon
Wants=network-online.target
After=network-online.target

[Service]
EnvironmentFile=-/etc/default/epoptes-client
ExecStartPre=/bin/sh -c "[ -s /etc/epoptes/server.crt ] || /usr/sbin/epoptes-client -c ${SERVER}"
ExecStart=/usr/sbin/epoptes-client

[Install]
WantedBy=multi-user.target
