# epoptes-client.service
[Unit]
Description=Epoptes-client system daemon
After=network.target

[Service]
ExecStart=/usr/sbin/epoptes-client

[Install]
WantedBy=multi-user.target
