# esekeyd.service
[Unit]
Description=ESE Key Daemon

[Service]
Type=forking
ExecStart=/usr/bin/esekeyd /etc/esekeyd.conf

[Install]
WantedBy=multi-user.target
