# etherpad-lite.socket
[Socket]
ListenStream=9000

[Install]
WantedBy=sockets.target
