# etherpad-lite.service
[Unit]
Description=Etherpad-lite
After=syslog.target network.target

[Service]
User=etherpad-lite
Group=etherpad-lite
StateDirectory=etherpad-lite
CacheDirectory=etherpad-lite
LogsDirectory=etherpad-lite
WorkingDirectory=/usr/share/etherpad-lite
ExecStart=/usr/bin/node /usr/share/etherpad-lite/src/node/server.js --settings /etc/etherpad-lite/settings.json
PrivateTmp=true
NoNewPrivileges=true
PrivateDevices=true
Restart=always

[Install]
WantedBy=multi-user.target
