# ezzram.service
[Service]
RemainAfterExit=yes
ExecStart=/usr/bin/ezzram start
ExecStop=/usr/bin/ezzram stop

[Install]
WantedBy=multi-user.target
