# fabric-peer.service
[Unit]
Description=Hyperledger Fabric Peer Node
After=network.target

[Service]
Type=simple
User=hyperledger
Group=hyperledger
StateDirectory=hyperledger hyperledger/db hyperledger/production
EnvironmentFile=/usr/lib/environment.d/fabric-peer.conf
ExecStart=/usr/bin/peer node start
PIDFile=/var/run/fabric-peer.pid

