# fabricd-backup.service
[Unit]
Description=Minecraft Fabric Server World Backup
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/fabricd backup
User=fabric
Group=fabric

[Install]
WantedBy=multi-user.target
