# fanctl.service
[Unit]
Description=fanctl fan controller

[Service]
Environment=RUST_LOG=fanctl=info
ExecStart=/usr/bin/fanctl -c /etc/fanctl.yml

[Install]
WantedBy=default.target
