# fancybeep.rules
# Give write access to the PC speaker only to the "input" group
KERNEL=="beep", ACTION=="add", SUBSYSTEM=="misc", GROUP="input", MODE="0640", RUN+="/usr/bin/fbdaemon -r"
                ACTION=="add", SUBSYSTEM=="input", ATTR{dev}=="13:*", ATTRS{name}=="Fancy Beeper", RUN+="/bin/sh -c '(sleep 1 && /bin/ln -sf $$DEVNAME /dev/input/by-path/platform-pcspkr-event-spkr) &'"
                ACTION=="remove", SUBSYSTEM=="input", ATTRS{name}=="Fancy Beeper", RUN+="/bin/rm -f /dev/input/by-path/platform-pcspkr-event-spkr"
