# fangfrisch.timer
[Unit]
Description=Download unofficial clamav virus definition files
Requires=fangfrisch.service

[Timer]
#OnCalendar=hourly
#RandomizedDelaySec=10m
OnCalendar=*:0/10
RandomizedDelaySec=2m
Unit=fangfrisch.service
Persistent=true

[Install]
WantedBy=timers.target
