# fangfrisch.timer
[Unit]
Description=Download unofficial clamav virus definition files
Requires=fangfrisch.service

[Timer]
OnCalendar=hourly
Unit=fangfrisch.service
RandomizedDelaySec=10m
Persistent=true

[Install]
WantedBy=timers.target
