# fanshim.service
[Unit]
Description=Fan Shim Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/fanshim --threshold 55 --hysteresis 5 --delay 2
Restart=on-failure

[Install]
WantedBy=multi-user.target
