# fanshim.service
[Unit]
Description=Fan Shim Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/fanshim --preempt --on-threshold 55 --off-threshold 50 --delay 2
Restart=on-failure

[Install]
WantedBy=multi-user.target
