# filebeat.service
[Unit]
Description=filebeat
Documentation=https://www.elastic.co/guide/en/beats/filebeat/current
After=network.target

[Service]
ExecStart=/usr/bin/filebeat -path.config /etc/filebeat \
                            -path.home /usr/share/filebeat \
                            -path.data /var/lib/filebeat \
                            -path.logs /var/log/filebeat
User=filebeat

[Install]
WantedBy=multi-user.target
