# filebrowser.service
[Unit]
Description=File Browser Service
After=network.target
Wants=network.target

[Service]
Type=simple
PIDFile=/var/run/filebrowser.pid
ExecStart=/usr/bin/filebrowser -c /etc/filebrowser/config.json
Restart=on-failure

[Install]
WantedBy=multi-user.target
