# firetable.service
[Unit]
Description=Firetable firewall daemon
After=network.target remote-fs.target nss-lookup.target

[Service]
Type=oneshot
ExecStart=/usr/bin/firetable start
ExecStop=/usr/bin/firetable stop
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
