# firetable.service
[Unit]
Description=Firetable firewall daemon

[Service]
Type=oneshot
ExecStart=/usr/bin/firetable start
ExecStop=/usr/bin/firetable stop
RemainAfterExit=true

[Install]
WantedBy=network.target
