# fix-bt-a2dp.service
[Unit]
Description=Fix Bluetooth A2DP

[Service]
Type=simple
ExecStart=/usr/bin/fix-bt-a2dp
StandardOutput=journal

[Install]
WantedBy=multi-user.target
