# flaresolverr.service
[Unit]
Description=FlareSolverr
After=network.target

[Service]
SyslogIdentifier=flaresolverr
Restart=always
RestartSec=5
Type=simple
User=flaresolverr
Group=flaresolverr
Environment="LOG_LEVEL=info"
Environment="CAPTCHA_SOLVER=none"
WorkingDirectory=/srv/http/flaresolverr
ExecStart=/usr/bin/python /srv/http/flaresolverr/src/flaresolverr.py
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
