# flexisip.service
[Unit]
Description=Flexisip server
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/flexisip --daemon -p /run/flexisip.pid
PIDFile=/run/flexisip.pid

[Install]
WantedBy=multi-user.target
