# flexisip.service
[Unit]
Description=Flexisip proxy server
After=syslog.target network.service network.target redis.service mysqld.service

[Service]
Type=forking
ExecStart=/usr/bin/flexisip --server proxy --daemon -p /run/flexisip.pid
PIDFile=/run/flexisip.pid

[Install]
WantedBy=multi-user.target
