# flood.service
[Unit]
Description=Flood rTorrent Web UI
After=network.target

[Service]
WorkingDirectory=/var/lib/flood
ExecStart=/usr/bin/npm start
User=flood

[Install]
WantedBy=multi-user.target
