# flow-causal.service
[Unit]
Description=flow causal service
After=syslog.target network.target

[Service]
ExecStart=/usr/bin/flow-causal
Type=simple


[Install]
WantedBy=default.target
