# foldingathome.service
[Unit]
Description=Folding@home distributed computing client
After=network.target

[Service]
Type=simple
User=fah
DynamicUser=yes
PermissionsStartOnly=true
ConfigurationDirectory=foldingathome
StateDirectory=fah
WorkingDirectory=/var/lib/fah
ExecStartPre=/usr/bin/chown fah:fah /etc/foldingathome/config.xml
ExecStart=/usr/bin/FAHClient --config /etc/foldingathome/config.xml
CPUSchedulingPolicy=idle
IOSchedulingClass=3

[Install]
WantedBy=multi-user.target

