# foldingathome-nvidia.service
[Unit]
Description=Folding@home helper for NVIDIA GPUs
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/bin/nvidia-modprobe
ExecStart=/usr/bin/nvidia-modprobe -c 0 -u
ExecStart=/usr/bin/clinfo
ExecStart=/usr/bin/sleep 60

[Install]
WantedBy=multi-user.target

