# foldingathome.service
[Unit]
Description=Folding@home distributed computing client
After=network.target

[Service]
Type=simple
User=fah
DynamicUser=yes
StateDirectory=fah
WorkingDirectory=/var/lib/fah
ExecStartPre=/usr/bin/cp /etc/foldingathome/config.xml /var/lib/fah/
ExecStart=/usr/bin/FAHClient
CPUSchedulingPolicy=idle
IOSchedulingClass=3

[Install]
WantedBy=multi-user.target

