# foldingathome.service
[Unit]
Description=Folding@home distributed computing client
After=network.target

[Service]
DynamicUser=yes
Type=simple
User=fah
SupplementaryGroups=video
StateDirectory=fah
WorkingDirectory=/var/lib/fah
ExecStart=/opt/fah/FAHClient

[Install]
WantedBy=multi-user.target
