# forged-backup.service
[Unit]
Description=Minecraft Forge 1.10.2 Server World Backup
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/bin/forge-1.10.2d backup
User=forge-1-10-2
Group=forge-1-10-2

[Install]
WantedBy=multi-user.target
