# forged.service
[Unit]
Description=Minecraft Forge 1.7.10 Server
After=local-fs.target network.target multi-user.target

[Service]
Type=forking
ExecStart=/usr/bin/forge-1.7.10d start
ExecStop=/usr/bin/forge-1.7.10d stop
User=forge-1-7-10
Group=forge-1-7-10

[Install]
WantedBy=multi-user.target
