# forged.service
[Unit]
Description=Minecraft Forge 1.8.9 Server
After=local-fs.target network.target multi-user.target

[Service]
Type=forking
ExecStart=/usr/bin/forge-1.8.9d start
ExecStop=/usr/bin/forge-1.8.9d stop
User=forge-1-8-9
Group=forge-1-8-9

[Install]
WantedBy=multi-user.target
