# forged-backup.service
[Unit]
Description=Minecraft Forge 1.9.4 Server World Backup
After=local-fs.target

[Service]
Type=forking
ExecStart=/usr/bin/forge-1.9.4d backup
User=forge-1-9-4
Group=forge-1-9-4

[Install]
WantedBy=multi-user.target
