# forged.service
[Unit]
Description=Minecraft Forge Server
After=local-fs.target network.target multi-user.target

[Service]
Type=forking
ExecStart=/usr/bin/forged start
ExecStop=/usr/bin/forged stop
User=forge
Group=forge

[Install]
WantedBy=multi-user.target
