# forti.service
[Unit]
Description=forticlient systemd

[Service]
ExecStart=/usr/lib/forti/forti.expect
EnvironmentFile=/etc/forti.conf
RemainAfterExit=true
