# fr24feed.service
[Unit]
Description=Flightradar24 decoder and feeder
After=network.target

[Service]
Restart=on-failure
User=fr24feed
ProtectSystem=full
ProtectHome=true
RuntimeDirectory=fr24feed

ExecStart=/usr/bin/fr24feed --config-file=/etc/fr24feed/fr24feed.ini --logfile=/var/log/fr24feed/fr24feed.log

[Install]
WantedBy=network.target

