# frpc.service
[Unit]
Description=Frp Client Service
After=network.target

[Service]
Type=exec
LoadCredential=config:/etc/frp/frpc.ini
DynamicUser=yes
Restart=on-failure
RestartSec=5s
ExecStart=/usr/bin/frpc -c ${CREDENTIALS_DIRECTORY}/config
ExecReload=/usr/bin/frpc reload -c ${CREDENTIALS_DIRECTORY}/config
# https://www.freedesktop.org/software/systemd/man/systemd.exec.html#Credentials

[Install]
WantedBy=multi-user.target
