# fs2es-indexer.service
[Unit]
Description=Filesystem indexer for Elasticsearch

[Service]
Type=simple
ExecStart=/usr/bin/fs2es-indexer daemon

[Install]
WantedBy=graphical.target
