# fsniper.service
[Unit]
Description=fsniper service

[Service]
ExecStart=/usr/bin/fsniper
Restart=always

[Install]
WantedBy=default.target
