# funkwhale.target
[Unit]
Description=Funkwhale
Wants=funkwhale-server.service funkwhale-worker.service funkwhale-beat.service

[Install]
WantedBy=multi-user.target
