# fusioninventory-agent.service
[Unit]
Description=FusionInventory Agent
After=network.service

[Service]
EnvironmentFile=/etc/conf.d/fusioninventory-agent
ExecStart=/usr/bin/fusioninventory-agent $FUSINVOPT

[Install]
WantedBy=multi-user.target
