# g13.service
[Unit]
Description=Logitech G13 daemon

[Service]
Type=simple
User=g13
Group=g13
ExecStart=/usr/bin/g13d --config /etc/g13/default.bind --pipe_in /run/g13d/g13-0 --pipe_out /run/g13d/g13-0_out &
