# gcp-cups-connector.service
[Unit]
Description=gcp-cups-connector startup script
Requires=org.cups.cupsd.service
After=network.target org.cups.cupsd.service

[Service]
Type=simple
ExecStart=/usr/bin/gcp-cups-connector --config-filename /etc/gcp-cups-connector/gcp-cups-connector.config.json

[Install]
WantedBy=multi-user.target
