# geoipupdate.timer
[Unit]
Description=Run GeoIP update weekly and on boot

[Timer]
OnBootSec=1min
OnUnitActiveSec=1w

[Install]
WantedBy=timers.target

