# gfwlist-route.timer
[Unit]
Description=GfwClient Update  Peer Half Hour  

[Timer]
Persistent=true
OnBootSec=2min
OnUnitActiveSec=30min
AccuracySec=1us
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target

