# gfwlist-route.service
[Unit]
Description=GFW Client Update
Requires=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/li/GitRepo/cmd/
ExecStart=/home/li/GitRepo/cmd/gfwclient
Restart=on-failure
#Restart=always
RestartSec=5

