# gfwlist-route.service
[Unit]
Description=GFW Client Update
Requires=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/lib/gfwlist-route/gfwlist-route
Restart=on-failure
#Restart=always
RestartSec=5

