# ghostfolio.service
[Unit]
Description=Open source wealth management
Wants=postgresql.service redis.service
After=postgresql.service redis.service

[Service]
Type=simple
ExecStart=/usr/bin/yarn start:prod
WorkingDirectory=/usr/lib/ghostfolio/api
EnvironmentFile=/etc/ghostfolio.env
User=ghostfolio
Group=ghostfolio
RestartSec=10

[Install]
WantedBy=multi-user.target
