# gitea.service
[Unit]
Description=Gogs (Go Git Service) server
After=syslog.target
After=network.target

[Service]
Type=simple
User=gogs
Group=gogs
ExecStart=/usr/share/gogs/gogs web --config /srv/gogs/conf/app.ini
WorkingDirectory=/usr/share/gogs/

[Install]
WantedBy=multi-user.target
