# github-actions.service
[Unit]
Description=GitHub Actions Listener Daemon

[Service]
User=github-actions
WorkingDirectory=~
Type=exec
ExecStart=./run.sh

[Install]
WantedBy=multi-user.target
