# github-actions.service
[Unit]
Description=GitHub Actions Listener Daemon

[Service]
User=github-actions
WorkingDirectory=/opt/github-actions
Type=exec
ExecStart=/opt/github-actions/run.sh

[Install]
WantedBy=multi-user.target
