# gitlab-ci-runner.service
[Unit]
Description=GitLab CI Runner

[Service]
User=gitlab_ci_runner
WorkingDirectory=/usr/lib/gitlab-ci-runner
SyslogIdentifier=gitlab-ci-runner
PIDFile=/run/gitlab-ci-runner.pid
ExecStart=/usr/bin/bundle exec "/usr/share/webapps/gitlab-ci-runner/bin/runner"
ExecStop=/usr/bin/kill -QUIT 
ExecReload=/usr/bin/kill -USR2 

[Install]
WantedBy=gitlab.target
