# gitlab-ci.target
[Unit]
Description=GitLab CI - Continuous integration server for gitlabhq
Requires=redis.service mysqld.service
After=redis.service mysqld.service syslog.target network.target gitlab.target

[Install]
WantedBy=multi-user.target
