# gitlab-ci-unicorn.service
[Unit]
Description=GitLab CI Unicorn Server

[Service]
User=gitlab_ci
WorkingDirectory=<DATADIR>
Environment=RAILS_ENV=production
SyslogIdentifier=gitlab-ci-unicorn
PIDFile=/run/gitlab-ci/unicorn.pid

ExecStart=/usr/bin/bundle exec "unicorn_rails -c <DATADIR>/config/unicorn.rb -E production"
ExecStop=/usr/bin/kill -QUIT $MAINPID
ExecReload=/usr/bin/kill -USR2 $MAINPID

[Install]
WantedBy=gitlab-ci.target
