# gksu-polkit.service
[Unit]
Description=Gksu server

[Service]
Type=dbus
BusName=org.gnome.Gksu
ExecStart=/usr/sbin/gksu-server
