# glance-registry.service
[Unit]
Description=glance Registry server

[Service]
ExecStart=/usr/bin/glance-registry

[Install]
WantedBy=multi-user.target
