# glassfish-web.service
[Unit]
Description = GlassFish Server v4 - Java EE 7 Web Profile
After = syslog.target network.target

[Service]
Type = forking
ExecStart = /opt/glassfish4-web/bin/asadmin start-domain domain1
ExecStop = /opt/glassfish4-web/bin/asadmin stop-domain domain1
ExecReload = /opt/glassfish4-web/bin/asadmin restart-domain domain1

[Install]
WantedBy = multi-user.target
