# globus-connect-personal.service
[Unit]
Description=Globus Connect Personal

[Service]
ExecStart=/usr/bin/globusconnect -start
ExecStop=/usr/bin/globusconnect -stop
Restart=on-failure

[Install]
WantedBy=default.target
