# glpi.service
[Unit]
Description=GLPI

[Service]
WorkingDirectory=/usr/share/webapps/glpi
ExecStart=/usr/bin/php -S localhost:8080 
