# glpi-cron.timer
[Unit]
Description=Timer for GLPI

[Timer]
OnUnitInactiveSec=1h

[Install]
WantedBy=timers.target
