# gmusicproxy.system.service
[Unit]
Description=Let's stream Google Play Music using any media-player

[Service]
ExecStart=/usr/bin/GMusicProxy -x -l /var/log/gmusicproxy.log -c /etc/gmusicproxy.cfg

[Install]
WantedBy=multi-user.target
