# gns3-server@.service
[Unit]
Description=GNS3 server for %i

[Service]
Type=simple
User=%i
ExecStart=/usr/bin/gns3server
Restart=on-failure

[Install]
WantedBy=multi-user.target
