# gnunet-user.service
[Unit]
Description=A framework for secure peer-to-peer networking

[Service]
Type=forking
ExecStart=gnunet-arm -c ~/.config/gnunet.conf -s
ExecStop=gnunet-arm -c ~/.config/gnunet.conf -e
GuessMainPID=no

[Install]
WantedBy=multi-user.target
