# gnunet.service
[Unit]
Description=A framework for secure peer-to-peer networking

[Service]
User=gnunet
Type=forking
ExecStart=/usr/bin/gnunet-arm -s -c /etc/gnunet.conf
ExecStop=/usr/bin/gnunet-arm -e -c /etc/gnunet.conf
StateDirectory=gnunet
GuessMainPID=no

[Install]
WantedBy=multi-user.target
