# go-hnrss.service
[Unit]
Description=Hacker News RSS
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/go-hnrss.conf
User=hnrss
ExecStart=/usr/bin/go-hnrss -bind ${LISTEN_ADDRESS} -rootRedirect ${ROOT_REDIRECT}
Restart=always

ProtectSystem=strict
PrivateDevices=true
ProtectKernelTunables=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target

