# go-nebulas.service
[Unit]
Description=Official Go implementation of the Nebulas protocol
Documentation=https://github.com/nebulasio/wiki
After=network.target

[Service]
ExecStart=/usr/bin/neb -c /opt/go-nebulas-git/conf/config.conf
WorkingDirectory=/opt/go-nebulas-git

[Install]
WantedBy=multi-user.target
