# go-wol-server.service
[Unit]
Description=Wake on Lan Service

[Service]
TimeoutStartSec=0
ExecStart=/usr/bin/go-wol-server -config=/etc/go-wol-server/config.toml

[Install]
WantedBy=multi-user.target