# go-wol-server.service
[Unit]
Description=Wake on Lan Service

[Service]
TimeoutStartSec=0
ExecStart=/usr/bin/go-wol-server -configFile=/etc/go-wol-server/config.toml -templateDir=/etc/go-wol-server/templates/ -remoteRoot=/ -listen=127.0.0.1:8080

[Install]
WantedBy=multi-user.target