# goaccess.timer
[Unit]
Description=Daily update of GoAccess Log Analysis

[Timer]
OnCalendar=hourly
RandomizedDelaySec=5min
Persistent=true

[Install]
WantedBy=timers.target
