# goaccess.service
[Unit]
Description=Goaccess Log Analysis
After=network.target

[Service]
Type=simple
User=root
Group=root
Restart=always
ExecStart=/bin/sh -c '/usr/bin/goaccess -f /var/log/nginx/*-access.log --enable-geopip --real-time-html -o /var/www/html/report.html'
[Install]
WantedBy=multi-user.target
