# gogs.service
[Unit]
Description=Gogs (Go Git Service) server
After=syslog.target
After=network.target

[Service]
Type=simple
User=gogs
Group=gogs
ExecStart=/usr/share/gogs/start.sh
WorkingDirectory=/usr/share/gogs/

[Install]
WantedBy=multi-user.target
