# goplay2.service
[Unit]
Description=Sound Service

Requires=pulseaudio.service
ConditionUser=!root

[Service]
ExecStart=/usr/bin/goplay2 -i enp5s0
Restart=on-failure

[Install]
Also=pulseaudio.service
WantedBy=default.target
