# goplay2.service
[Unit]
Description=Sound Service

Requires=pulseaudio.service
ConditionUser=!root

[Service]
ExecStart=/usr/bin/goplay2 -i enp5s0 -sink alsa_output.pci-0000_43_00.1.hdmi-stereo-extra3
Restart=on-failure

[Install]
Also=pulseaudio.service
WantedBy=default.target
