# goplay2.service
[Unit]
Description=Sound Service

Requires=pulseaudio.service
ConditionUser=!root

[Service]
ExecStart=/usr/bin/goplay2 -i $(cat $HOME/.config/goplay2-network-interface || cat /etc/goplay2-default-network-interface)
Restart=on-failure

[Install]
Also=pulseaudio.service
WantedBy=default.target
