# gopls.service
[Unit]
Description=Language Server Protocol implementation for Go
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/gopls

[Install]
WantedBy=multi-user.target
