# gosearch.install
post_install() {
  systemctl enable gosearch.service
  systemctl start gosearch.service
}

pre_upgrade() {
  systemctl stop gosearch.service
}

post_upgrade() {
  systemctl daemon-reload
  systemctl start gosearch.service
}

pre_remove() {
  systemctl stop gosearch.service
}
