# gotify-server.service
[Unit]
Description=Gotify Server
After=network.target

[Service]
Type=simple
User=gotify
Group=gotify
ExecStart=/usr/bin/gotify-server

[Install]
WantedBy=multi-user.target
