# gpgsync.timer
[Unit]
Description=Periodically sync GPG keys

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=2000

[Install]
WantedBy=timers.target
