# gpmdp2text.service
[Unit]
Description=Utility to read GPMDP's JSON API and output as plain text

[Service]
Type=simple
ExecStart=/usr/bin/gpmdp2text

[Install]
WantedBy=default.target