# graphite-web.service
[Unit]
Description=Graphite-web daemon
After=network.target

[Service]
Environment=PYTHONPATH=/opt/graphite/whisper
ExecStart=/opt/graphite/bin/run-graphite-devel-server.py --libs=/opt/graphite/webapp/ /opt/graphite/
Type=simple

# some hardening
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true
CapabilityBoundingSet=
ReadOnlyPaths=/

[Install]
WantedBy=multi-user.target

