# graphite-exporter.service
[Unit]
Description=Prometheus exporter for Graphite metrics
Requires=network-online.target
After=network-online.target

[Service]
EnvironmentFile=-/etc/conf.d/graphite-exporter
ExecStart=/usr/bin/graphite_exporter $ARGS
ExecReload=/bin/kill -HUP $MAINPID
DynamicUser=yes
Restart=on-failure
NoNewPrivileges=true
ProtectHome=read-only
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
