# graviumd@.service
[Unit]
Description=Gravium Daemon
After=network.target

[Service]
User=%i
Group=users
Type=simple
ExecStart=/usr/bin/graviumd

[Install]
WantedBy=multi-user.target
