# greptimedb.service
[Unit]
Description=GreptimeDB time series database
Documentation=https://docs.greptime.com/
After=network.target

[Service]
Type=simple
User=greptime
ExecStart=/usr/bin/greptime --log-dir /var/lib/greptimedb/logs standalone start -c /etc/greptimedb/greptimedb.toml
ExecStop=/usr/bin/greptime standalone stop

[Install]
WantedBy=multi-user.target
