# groovebasin.service
[Unit]
Description=groovebasin music player
Documentation=http://groovebasin.com/

[Service]
User=groovebasin
ExecStart=/usr/bin/groovebasin
WorkingDirectory=/var/lib/groovebasin

[Install]
WantedBy=multi-user.target
