# guacd.service
[Unit]
Description=Guacamole Proxy Daemon
After=network.target

[Service]
Type=forking
PIDFile=/var/run/guacd.pid
ExecStart=/usr/bin/guacd -p /var/run/guacd.pid
Restart=on-failure

[Install]
WantedBy=multi-user.target

