# hasplmd.service
[Unit]
Description=AKSUSB winehasp daemon
After=aksusbd.service winehasp.service
Requires=aksusbd.service winehasp.service
BindsTo=aksusbd.service winehasp.service
PartOf=aksusbd.service winehasp.service

[Service]
Type=forking
StandardOutput=tty
TimeoutStopSec=10s

PrivateTmp=true
PrivateNetwork=true
ProtectSystem=full
ProtectHome=true

ExecStart=/usr/bin/hasplmd -s

[Install]
WantedBy=multi-user.target
Also=aksusbd.service
Also=winehasp.service
