# headphones-system.service
[Unit]
Description=Headphones Daemon
After=network.target

[Service]
User=headphones
Type=forking
PIDFile=/run/headphones/headphones.pid
ExecStart=/usr/bin/env python2 /opt/headphones/Headphones.py -q --daemon --nolaunch --pidfile=/run/headphones/headphones.pid --config /opt/headphones/config.ini --datadir /opt/headphones

[Install]
WantedBy=multi-user.target
