# headscale.service
[Unit]
Description=Headscale server
After=network.target

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/bin/headscale serve
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target
