# headphones-user.service
[Unit]
Description=Headphones Daemon

[Service]
ExecStart=/usr/bin/python2 /usr/lib/headphones/Headphones.py --config %h/.headphones/config.ini --datadir %h/.headphones

[Install]
WantedBy=default.target
