# headphones-user.service
[Unit]
Description=Headphones Daemon

[Service]
ExecStart=/usr/bin/env python2 /opt/headphones/Headphones.py --quiet --config %h/.headphones/config.ini --datadir %h/.headphones

[Install]
WantedBy=default.target
