# hentaiathome@.service
[Unit]
Description=Hentai At Home for User %i
After=network-online.target

[Service]
Type=simple
User=%i
ExecStart=/usr/bin/hentaiathome-cli

[Install]
WantedBy=multi-user.target
