# hexagonrpcd-adsp-rootpd.service
[Unit]
Name=Hexagonrpcd ADSP RootPD
Description=Daemon to support Qualcomm Hexagon ADSP virtual filesystem for RootPD
Requires=dev-fastrpc\x2dadsp.device
After=dev-fastrpc\x2dadsp.device

[Service]
User=fastrpc
Group=fastrpc
ExecStart=/usr/bin/hexagonrpcd -f /dev/fastrpc-adsp

[Install]
WantedBy=multi-user.target
