# hexagonrpcd-adsp-sensorspd.service
[Unit]
Description="Hexagonrpcd ADSP SensorPD"

[Service]
User=fastrpc
Group=fastrpc
ExecStart=/usr/bin/hexagonrpcd -f /dev/fastrpc-adsp -s

[Install]
WantedBy=multi-user.target
