# hexagonrpcd-sdsp.service
[Unit]
Name=Hexagonrpcd SDSP
Description=Daemon to support Qualcomm Hexagon SDSP virtual filesystem
Requires=dev-fastrpc\x2dsdsp.device
After=dev-fastrpc\x2dsdsp.device

[Service]
User=fastrpc
Group=fastrpc
ExecStart=/usr/bin/hexagonrpcd -f /dev/fastrpc-sdsp -s

[Install]
WantedBy=multi-user.target
