# hexagonrpcd-adsp-sensorspd.service
[Unit]
Name=Hexagonrpcd ADSP SensorPD
Description=Daemon to support Qualcomm Hexagon ADSP virtual filesystem for SensorPD
Requires=dev-fastrpc\x2dadsp.device
After=dev-fastrpc\x2dadsp.device

[Service]
User=fastrpc
Group=fastrpc
ExecStart=/usr/bin/hexagonrpcd -f /dev/fastrpc-adsp -s

[Install]
WantedBy=multi-user.target
