# hhvm@.service
[Unit]
Description=HHVM FastCGI application server
After=syslog.target network.target

[Service]
EnvironmentFile=/etc/hhvm/systemd-%I.conf

PIDFile=/run/hhvm/hhvm-%I.pid
PrivateTmp=true
ExecStart=/usr/bin/hhvm --mode server --config=/etc/hhvm/%I.ini --user ${User} -vPidFile=/run/hhvm/hhvm-%I.pid

[Install]
WantedBy=multi-user.target
