# hiawatha-letsencrypt-renew.service
[Unit]
Description=Renew Let's Encrypt certificates for Hiawatha
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/share/hiawatha-letsencrypt/letsencrypt renew restart
