# himitsud.service
[Unit]
Description=Himitsu Daemon

[Service]
ExecStart=/usr/bin/himitsud
Restart=always

[Install]
WantedBy=default.target
