# hissh-agent.service
[Unit]
Description=Himitsu SSH Agent

[Service]
ExecStart=/usr/bin/hissh-agent
Restart=always

[Install]
WantedBy=default.target
