# homegear.service
[Unit]
Description=Program to interface smart home devices with home automation software
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/homegear
ExecReload=/usr/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
