# homebridge-system.service
[Unit]
Description=Node.js HomeKit Server
Requires=avahi-daemon.service
After=network.target avahi-daemon.service

[Service]
User=homebridge
Type=simple
ExecStart=/usr/bin/homebridge -U /var/lib/homebridge

[Install]
WantedBy=multi-user.target
