# homeassistant-printerqueue@.service
[Unit]
Description=Turn printer automatically on when printing using homeassistant

[Service]
EnvironmentFile=/etc/homeassistant/%i.conf
ExecStart=/usr/lib/systemd/system/homeassistant-printerqueue.sh
Restart=always

[Install]
WantedBy=multi-user.target
