# honk.service
[Unit]
Description="Honk ActivityPub server"

[Service]
WorkingDirectory=/usr/share/honk
ExecStart=/usr/bin/honk

User=honk
Group=honk

[Install]
WantedBy=default.target
