# hopglass-server@.service
[Unit]
Description=The HopGlass Server collects data from the Freifunk network (%I)

[Service]
Type=simple
User=hopglass-server
Group=nobody
StateDirectory=hopglass-server/%I
WorkingDirectory=/var/lib/hopglass-server/%I/
ExecStart=/usr/bin/hopglass-server --config /etc/hopglass-server/%I/config.json
Restart=always

[Install]
WantedBy=multi-user.target
