# hornet.service
[Unit]
Description=Hornet
Wants=network-online.target
After=network-online.target

[Service]
LimitNOFILE=4096
SyslogIdentifier=hornet
PrivateDevices=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes

User=hornet
WorkingDirectory=/var/lib/hornet
TimeoutSec=1200
Restart=always
ExecStart=/usr/bin/hornet --config "/etc/hornet/config.json" --peering "/etc/hornet/peering.json" --profiles "/etc/hornet/profiles.json"

[Install]
WantedBy=multi-user.target
