# hosts-gen.path
[Unit]
Description=Automatically regenerate hosts from hosts.d
Documentation=file:/usr/share/doc/hosts-gen/README

[Path]
PathChanged=/etc/hosts.d
Unit=hosts-gen.service

[Install]
WantedBy=default.target
