# hosts.service
[Service]
Type=oneshot
ExecStart=/usr/bin/curl -s --retry 1024 -o /etc/hosts "https://raw.githubusercontent.com/StevenBlack/hosts/master/hosts"
