# hottub.service
[Unit]
Description=CI bridge from GitHub to SourceHut
After=network.target

[Service]
Type=simple
User=hottub
Restart=always
EnvironmentFile=/etc/hottub/params
ExecStart=/usr/bin/hottub -listen=${listen} -db=/var/lib/hottub/db -gh-app-id=${gh_app_id} -gh-private-key=${gh_private_key} -gh-webhook-secret=${gh_webhook_secret} -metasrht-client-id=${metasrht_client_id} -metasrht-client-secret=${metasrht_client_secret}

[Install]
WantedBy=multi-user.target
