# hqplayerd.service
[Unit]
Description=HQPlayer Embedded daemon
Requires=network-online.target sound.target
After=network-online.target sound.target local-fs.target remote-fs.target

[Service]
Type=simple
ExecStartPre=/bin/sleep 60
ExecStart=/usr/bin/hqplayerd
Restart=always
RestartSec=2
EnvironmentFile=-/etc/default/hqplayerd

[Install]
WantedBy=multi-user.target
