# hqplayerd_rygel.service
[Unit]
Description=HQPlayer Embedded daemon with rygel
After=network.target sound.target

[Service]
Type=simple
Restart=always
ExecStart=/bin/bash /usr/bin/start_hqplayerd.sh

[Install]
WantedBy=multi-user.target