# hqplayerd_rygel.service
[Unit]
Description=HQplayer embedded
After=network.target sound.target

[Service]
ExecStart=/bin/bash /usr/bin/start_hqplayerd.sh

[Install]
WantedBy=default.target