# hqplayerd.service
[Unit]
Description=HQPlayer Embedded daemon
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/hqplayerd
Restart=always
#EnvironmentFile=-/etc/default/hqplayerd

[Install]
WantedBy=default.target
