# httpjs.service
[Unit]
Description=A web server written in NodeJS
After=network.target

[Service]
Type=forking
PIDFile=/run/httpjs.pid
PrivateDevices=yes
SyslogLevel=err

ExecStart=/home/py64/GIT/HTTPjs/http.js --pid /run/httpjs.pid
ExecReload=/usr/bin/kill -HUP $MAINPID
KillSignal=SIGQUIT
KillMode=mixed

[Install]
WantedBy=multi-user.target

