# https_dns_proxy.service
[Unit]
Description=Light-weight DNS<-->HTTPS, non-caching proxy
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
AmbientCapabilities=CAP_NET_BIND_SERVICE
NoNewPrivileges=true
ExecStart=/usr/bin/https_dns_proxy -u nobody -g nobody -v
Restart=on-failure

[Install]
WantedBy=multi-user.target
