# hushd.service
[Unit]
Description=Hush daemon
After=network.target

[Service]
ExecStart=/usr/bin/hushd

[Install]
WantedBy=default.target
