# hylafax.cron.daily
#!/bin/sh
# HylaFAX+

/usr/bin/faxqclean
/usr/bin/faxcron -rcv 30
