# hypervfcopyd.service
[Unit]
Description=Hyper-V FIle Copy daemon
ConditionalVirtualization=microsoft

[Service]
Type=simple
ExecStart=/usr/local/bin/hypervfcopyd -n

[Install]
WantedBy=multi-user.target
