# hyperv-enhanced-session.service
[Unit]
Description=Hyper-V Enhanced Session Mode Daemon
After=network-online.target
Before=xrdp.service xrdp-sesman.service
Wants=network-online.target xrdp.service xrdp-sesman.service

[Service]
Type=oneshot
ExecStart=/opt/%N/configure-rdp.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
