# hyperiond.service
[Unit]
Description=HyperionD
Before=vdr.service

[Service]
ExecStart=/usr/bin/hyperiond /etc/hyperion.config.json
Restart=on-abort

[Install]
WantedBy=multi-user.target