# hypervvssd.service
[Unit]
Description=Hyper-V VSS daemon
ConditionalVirtualization=microsoft

[Service]
Type=simple
ExecStart=/usr/local/bin/hypervvssd -n

[Install]
WantedBy=multi-user.target
