# hysteria1-server@.service
[Unit]
Description=Hysteria (Legacy) Server Service
Documentation=https://v1.hysteria.network/
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/hysteria1 server --config %i.json --no-check
WorkingDirectory=/etc/hysteria
User=hysteria
Group=hysteria
Environment=HYSTERIA_LOG_LEVEL=info
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_BIND_SERVICE CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_BIND_SERVICE CAP_NET_RAW
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
