# i2pd.service
[Unit]
Description=I2P Router in C++
After=network.target

[Service]
User=i2p
Group=i2p
Type=forking
ExecStart=/usr/bin/i2p --conf=/var/lib/i2pd/i2p.conf --service=1 --daemon=1
ExecReload=/usr/bin/kill -HUP $MAINPID
PIDFile=/run/i2pd.pid
KillSignal=SIGINT
LimitNOFILE=8192
PrivateDevices=yes

[Install]
WantedBy=multi-user.target
