# i2pd.service
[Unit]
Description=C++ daemon for accessing the I2P network
After=network.target

[Service]
Type=forking
Restart=on-abnormal
PIDFile=/run/i2pd/i2pd.pid
User=i2pd
Group=i2pd
LimitNOFILE=4096
PermissionsStartOnly=yes
ExecStart=/usr/bin/i2pd --daemon --service --pidfile=/run/i2pd/i2pd.pid --conf=/etc/i2pd/i2pd.conf --tunconf=/etc/i2pd/tunnels.conf

[Install]
WantedBy=multi-user.target
