# i8kmon.service
[Unit]
Description=Dell laptop thermal monitoring
ConditionPathExists=/proc/i8k

[Service]
ExecStartPre=/usr/bin/smm 30a3  # uncomment to disable BIOS fan control
ExecStopPost=/usr/bin/smm 31a3  # ... and re-enable it afterwards
ExecStart=/usr/bin/i8kmon -d

[Install]
WantedBy=multi-user.target
