# ibus-config@.service
[Unit]
Description=Intelligent Input Bus Config Daemon

[Service]
Slice=ibus.slice
ExecStartPre=/bin/sleep 2
ExecStart=/usr/lib/ibus/ibus-dconf
Environment=DISPLAY="%I" GTK_IM_MODULE=ibus QT_IM_MODULE=ibus XMODIFIERS="@im=ibus"
