# ibus@.service
[Unit]
Description=Intelligent Input Bus
Documentation=man:ibus-daemon(1)
Conflicts=fcitx.service

[Service]
ExecStart=/usr/bin/ibus-daemon --replace --xim
ExecReload=/usr/bin/ibus restart
ExecStop=/usr/bin/ibus exit
Environment=DISPLAY="%I" GTK_IM_MODULE=ibus QT_IM_MODULE=ibus XMODIFIERS="@im=ibus"

[Install]
Alias=input-method.service
WantedBy=default.target
