# icecream-scheduler.service
[Unit]
Description=Icecream distributed compiler scheduler

[Service]
Type=simple
User=icecream
Group=icecream
SyslogIdentifier=icecc-scheduler
ExecStart=/usr/lib/icecream/icecream-schedulerd

[Install]
WantedBy=multi-user.target
