# iked.service
[Unit]
Description=OpenIKED
Requires=network-online.target
After=network-online.target

[Service]
Type=forking
ExecStart=/usr/local/sbin/iked

[Install]
WantedBy=multi-user.target
