# iked.service
[Unit]
Description=iked daemon
After=syslog.target network.target

[Service]
Type=forking
ExecStart=/usr/sbin/iked

[Install]
WantedBy=multi-user.target
