# indy-node-control.service
[Unit]
Description=Service for upgrade of existing Indy Node and other operations
#Requires=indy.service
#After=indy.service
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/indy/node_control.conf
ExecStart=/usr/bin/start_node_control_tool $TEST_MODE --hold-ext $HOLD_EXT
Restart=on-failure
RestartSec=10
StartLimitBurst=10
StartLimitInterval=200
TimeoutSec=300

[Install]
WantedBy=multi-user.target
