# initrd-network.service

# NOTE: to be included in initramfs the file must contain "/etc/initrd-release" string

[Unit]
Description=Network Service (initrd)
ConditionPathExists=/etc/initrd-release
DefaultDependencies=no
Requires=systemd-networkd.service
Requires=systemd-resolved.service
After=network.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/bin/true

[Install]
WantedBy=sysinit.target
