# initrd-dropbear.service

# NOTE: to be included in initramfs the file must contain "/etc/initrd-release" string

[Unit]
Description=Dropbear SSH Server (initrd)
ConditionPathExists=/etc/initrd-release
#DefaultDependencies=no
Requires=initrd-network.service
After=network.target

[Service]
ExecStart=/usr/bin/dropbear -E -s -j -k

[Install]
WantedBy=sysinit.target
