# invoiceninja-scheduler.timer
[Unit]
Description=Invoice Ninja scheduler timer

[Timer]
# every minute
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
