# ipfs-fuse.service
[Unit]
Description=Daemon to mount IPFS
Conflicts=ipfs.service

[Service]
ExecStart=/usr/bin/start-ipfs-fuse-daemon.sh
Restart=on-failure

[Install]
WantedBy=default.target
