# ipfs.service
[Unit]
Description=IPFS daemon
After=network-online.target
Requires=network-online.target

[Service]
User=ipfs
Environment="IPFS_PATH=/opt/ipfs-repo"
ExecStart=/usr/bin/ipfs daemon

[Install]
WantedBy=multi-user.target
