# irexec.service
[Unit]
Description=irexec as service

[Service]
ExecStart=/usr/bin/irexec
Restart=on-abort

[Install]
WantedBy=default.target

