# isso.service
[Unit]
Description=Isso Comment Server

[Service]
Environment="ISSO_SETTINGS=/etc/isso.conf"
ExecStart=/usr/bin/isso run

[Install]
WantedBy=multi-user.target
