# isso.service
[Unit]
Description=Isso Comment Server

[Service]
EnvironmentFile=/etc/conf.d/isso
ExecStart=/usr/bin/isso -c $ISSO_CONF

[Install]
WantedBy=multi-user.target

