# isso.service
[Unit]
Description=Isso Commenting Server
After=network.target network-online.target nss-lookup.target

[Service]
User=isso
Group=isso
ExecStart=/usr/bin/isso -c /etc/isso.conf run

[Install]
WantedBy=multi-user.target
