# iwatch.service
[Unit]
Description=iwatch daemon

[Service]
ExecStart=/usr/bin/iwatch -f /etc/iwatch/iwatch.xml

[Install]
WantedBy=multi-user.target
