# iyuuplus.service
[Unit]
Description = IYUU Auto Reseed Plus Daemon.
Wants = network-online.target php-fpm.service
After = network-online.target php-fpm.service

[Service]
User = iyuuplus
Group = iyuuplus
Restart = on-failure
ExecStart = /usr/bin/iyuuplus start
ExecStop = /usr/bin/iyuuplus stop -g
ExecReload = /usr/bin/iyuuplus reload -g

[Install]
WantedBy = multi-user.target
