# jackd@.service
[Unit]
Description=JACK low-latency audio server daemon
After=sound.target

[Service]
EnvironmentFile=/etc/conf.d/jackd.conf
User=%I
ExecStart=/usr/bin/jackd $OPTIONS

[Install]
WantedBy=multi-user.target
