# jdownloader.service
[Unit]
Description=JDownload Headless
After=network.target

[Service]
ExecStart=/usr/bin/JDownloaderHeadless
ExecStop=
User=jdownloader
Group=jdownloader
Restart=on-success
[Install]
WantedBy=multi-user.target
