# jenkins.service
[Unit]
Description=Extendable continuous integration server
After=network.target

[Service]
User=jenkins
Type=exec
EnvironmentFile=/etc/conf.d/jenkins
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=jenkins
ExecStart=/bin/sh -c 'eval $JENKINS_COMMAND_LINE'

[Install]
WantedBy=multi-user.target
