# jetty.service
# Systemd unit file for jetty
#
# Multiple copies of this service (i.e. multiple concurrently running
# jetty servers) are not supported right now. Expect this to come with
# future updates

[Unit]
Description=Jetty Web Application Server
After=syslog.target network.target

[Service]
Type=forking
ExecStart=/usr/bin/jetty start
ExecStop=/usr/bin/jetty stop
ExecReload=/usr/bin/jetty restart
PIDFile=/run/jetty/jetty.pid

[Install]
WantedBy=multi-user.target
