# juiced.service
[Unit]
Description=Juiced battery daemon
Documentation=https://github.com/drkhsh/juiced

[Service]
Type=simple
ExecStart=/usr/bin/juiced
Restart=on-failure

[Install]
WantedBy=default.target
