# kairosdb.service
[Unit]
Description=KairosDB is a fast distributed scalable time series database written on top of Cassandra.
Requires=cassandra.service
After=cassandra.service

[Service]
Type=forking
User=root
PIDFile=/var/run/kairosdb.pid
ExecStart=/opt/kairosdb/bin/kairosdb.sh start
ExecStop=/opt/kairosdb/bin/kairosdb.sh stop

[Install]
WantedBy=multi-user.target