# kamailio.service
[Unit]
Description=Kamailio - the Open Source SIP Server
Wants=network-online.target
After=syslog.target network-online.target

[Service]
Type=forking
PIDFile=/run/kamailio.pid
ExecStart=/usr/bin/kamailio -P /run/kamailio.pid
KillMode=mixed

[Install]
WantedBy=multi-user.target
