# kanboard.service
[Unit]
Description=Run kanboard's cron job

[Service]
Type=oneshot
WorkingDir=/usr/share/webapps/kanboard
ExecStart=/usr/share/webapps/kanboard/cli cronjob
