# kbct.service
[Unit]
Description=Keyboard keycode mapping daemon supporting layered configuration

[Service]
Type=simple
ExecStart=kbct remap --config /etc/kbct/config.yml
Restart=always

[Install]
WantedBy=default.target
