# kbpwrd.service
[Unit]
Description PinePhone (Pro) Keyboard Case Power Manager
[Service]
Type=simple
ExecStart=/usr/bin/kbpwrd
[Install]
WantedBy=sysinit.target

