# kcrap.service
[Unit]
Description=Kerberos Challenge-Response Authentication Protocol (KCRAP)
DefaultDependencies=no
After=krb5-kdc.service

[Service]
Type=forking
PIDFile=/run/kcrap_server.pid
ExecStart=/usr/sbin/kcrap_server
ExecStop=/bin/kill $MAINPID

[Install]
WantedBy=multi-user.target
