# kd-server.service
[Unit]
Description=kd the command-line dictionary's server

[Service]
Type=simple
ExecStart=/usr/bin/kd --server
Restart=always

[Install]
WantedBy=default.target

