# keydb-sentinel.service
[Unit]
Description=High availability daemon for KeyDB
After=network.target

[Service]
Type=notify
User=keydb
Group=keydb
ExecStart=/usr/bin/keydb-sentinel /etc/keydb/sentinel.conf --supervised systemd
TimeoutStartSec=60
TimeoutStopSec=60
CapabilityBoundingSet=
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/etc/keydb/
NoNewPrivileges=true
RuntimeDirectory=keydb
RuntimeDirectoryMode=755
LimitNOFILE=10032

[Install]
WantedBy=multi-user.target
