# keymapperd.service
[Unit]
Description=Keymapper Daemon

[Service]
ExecStart=/usr/bin/keymapperd

[Install]
WantedBy=multi-user.target
