# kibana.service
[Unit]
Description=Kibana - dashboard for Elasticsearch
After=elasticsearch.service

[Service]
User=nobody
Group=elasticsearch
Environment="NODE_ENV=production" "BABEL_CACHE_PATH=/tmp/kibana-babel-cache"
ExecStart=/usr/bin/node /usr/lib/kibana/src/cli --quiet --config=/etc/elasticsearch/kibana/kibana.yml

[Install]
WantedBy=multi-user.target
