# kibana.service
[Unit]
Description=Kibana - dashboard for Elasticsearch
After=elasticsearch.service

[Service]
User=nobody
Environment="CONFIG_PATH=/etc/elasticsearch/kibana/kibana.yml" "NODE_ENV=production"
ExecStart=/usr/bin/node /usr/lib/kibana/bin/kibana.js --quiet

[Install]
WantedBy=multi-user.target
