# klipper.service
[Unit]
Description=3D printer firmware with motion planning on the host

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
User=klipper
RemainAfterExit=yes
ExecStart=/usr/bin/python2 /opt/klipper/klippy/klippy.py /etc/klipper/klipper.cfg -l /var/log/klipper/klippy.log -I /run/klipper/sock
Restart=always
RestartSec=10
