# kmodcache.service
[Unit]
Description=update modules cache for active kernel
RequiresMountsFor=/var/cache

[Service]
Type=oneshot
ExecStart=/usr/bin/make-kmod-cache

[Install]
WantedBy=multi-user.target
