# kolide-launcher.service
[Unit]
Description=The Kolide Launcher
After=network.service syslog.service

[Service]
ExecStart=/usr/bin/kolide-launcher \
          -config /etc/kolide/launcher.conf
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
