# ktproxy-client@.service
[Unit]
Description=ktproxy service
After=network.target

[Service]
ExecStart=/usr/bin/ktproxy-client -c /etc/ktproxy/%i.toml
MemoryAccounting=true

[Install]
WantedBy=multi-user.target

