# kurtosis.service
[Unit]
Description=Starts kurtosis service

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/kurtosis engine start
ExecReload=/usr/bin/kurtosis engine restart
KillMode=process
TimeoutSec=180

[Install]
WantedBy=default.target
