# lanraragi.service
[Unit]
Description=LANraragi Service
After=syslog.target network.target redis.service
Requires=redis.service

[Service]
User=lanraragi
Group=lanraragi
Environment=PERL5LIB="/usr/lib/lanraragi/perl5/lib/perl5:$PERL5LIB"
WorkingDirectory=/usr/lib/lanraragi
ExecStart=/usr/bin/perl /usr/lib/lanraragi/script/launcher.pl -f /usr/lib/lanraragi/script/lanraragi
Type=simple
TimeoutStopSec=20
SyslogIdentifier=lanraragi
KillMode=control-group

[Install]
WantedBy=multi-user.target
