# lantern.service
[Unit]
Description=Lantern
After=network.target

[Service]
ExecStart=/usr/bin/lantern

[Install]
WantedBy=default.target
