# lazylibrarian.service
[Unit]
Description=LazyLibrarian

[Service]
ExecStart=/usr/bin/python2 /usr/lib/lazylibrarian/LazyLibrarian.py --config /etc/lazylibrarian/lazylibrarian.ini --datadir /var/lib/lazylibrarian --nolaunch
Type=simple
User=lazylibrarian
Group=lazylibrarian
Restart=always

[Install]
WantedBy=multi-user.target
