# leapd.service
[Unit]
Description=Leap Motion Service Daemon

[Service]
PIDFile=/var/run/leapd.pid
ExecStart=/usr/bin/leapd
KillSignal=SIGKILL
TimeoutStopSec=5
KillMode=process

[Install]
WantedBy=multi-user.target
