# leapd.service
[Unit]
Description=Leap Motion Service Daemon

[Service]
ExecStart=/usr/bin/leapd

[Install]
WantedBy=multi-user.target
