# lesspassgo-server.service
[Unit]
Description=LessPass API compatible server written in Go
Documentation=https://github.com/tuxlinuxien/lesspassgo/blob/master/README.md
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
DynamicUser=yes
StateDirectory=lesspassgo-server
EnvironmentFile=/etc/lesspassgo-server.conf
ExecStart=/usr/bin/lesspassgo-server

[Install]
WantedBy=multi-user.target
