# libinput-xrandr-autorotate.service
[Unit]
Description=Allows to rotate screen and touch input automatically

[Service]
Type=idle
ExecStart=/usr/bin/libinput-xrandr-autorotate
Restart=always

[Install]
WantedBy=multi-user.target
