# libinput-xrandr-autorotate.service
[Unit]
Description=libinput-xrandr-autorotate

[Service]
ExecStart=/usr/bin/libinput-xrandr-autorotate
Restart=always

[Install]
WantedBy=multi-user.target
