# libresonic.service
[Unit]
Description=Libresonic
After=network.target
 
[Service]
User=libresonic
ExecStart=/var/lib/libresonic/libresonic.sh
ExecStop=/bin/kill -15 $MAINPID
Type=forking
 
[Install]
WantedBy=multi-user.target
