# libv4l-rkmpp-setup.service
[Unit]
Description=Setup dummy V4L2 devices for chromium MPP decoding/encoding

[Service]
Type=oneshot
EnvironmentFile=/etc/conf.d/libv4l-rkmpp
ExecStart=/usr/bin/libv4l-rkmpp-setup.sh $PROFILE

[Install]
WantedBy=multi-user.target