# libvirt-snmp.service
[Unit]
Description=libvirt SNMP sub-agent
Requires=snmpd.service
After=snmpd.service

[Service]
Type=simple
ExecStart=/usr/bin/libvirtMib_subagent -f -L

[Install]
WantedBy=multi-user.target
