# lidarr.service
[Unit]
Description=Lidarr Service
After=syslog.target network.target

[Service]
User=lidarr
Group=lidarr
ExecStart=/usr/lib/lidarr/Lidarr -nobrowser -data=/var/lib/lidarr
Type=simple
TimeoutStopSec=20

# Lets built in updater work well.
Restart=on-failure
KillMode=process

[Install]
WantedBy=multi-user.target
