# lidarr.service
[Unit]
Description=Lidarr Service
After=syslog.target network.target

[Service]
User=lidarr
Group=lidarr
ExecStart=/usr/lib/lidarr/bin/Lidarr -nobrowser -data=/var/lib/lidarr
Type=simple

[Install]
WantedBy=multi-user.target
