# lispd.service
[Unit]
Description=LISPmob service

[Service]
ExecStart=/usr/bin/lispd -f /etc/lispd.conf

[Install]
WantedBy=multi-user.target

