# lispd.service
[Unit]
Description=LISPmob service

[Service]
ExecStart=/usr/bin/lispd

[Install]
WantedBy=multi-user.target
