# lldpd.service
[Unit]
Description=LLDP daemon
Documentation=man:lldpd(8)
Requires=snmpd.service
After=snmpd.service network.target
RequiresMountsFor=/run

[Service]
Type=notify
NotifyAccess=main
ExecStart=/usr/bin/lldpd -x -c -f -s -e -l -M 2
Restart=on-failure

[Install]
WantedBy=multi-user.target
