# llmnrd.service
[Unit]
Description=Link-Local Multicast Resolution Daemon
Before=systemd-resolved.service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/llmnrd -6
Restart=on-failure

[Install]
WantedBy=multi-user.target
