# local-ai.service
[Unit]
Description=LocalAI API For Running LLaMA/GPT Models
After=local-fs.target

[Service]
User=local-ai
Group=local-ai
EnvironmentFile=/etc/local-ai/local-ai.conf
SyslogIdentifier=local-ai
Type=simple
ExecStart=/usr/bin/local-ai
SuccessExitStatus=1

[Install]
WantedBy=multi-user.target
