# lockd.service
[Unit]
Description=Lockscreen daemon

[Service]
ExecStart=/usr/bin/lockd

[Install]
WantedBy=default.target