# logstash.service
[Unit]
Description=logstash
After=elasticsearch.service

[Service]
Type=simple
User=logstash
Group=logstash
Environment=LOGSTASH_CONF=/etc/logstash/
ExecStart=/usr/share/logstash/bin/logstash "--path.settings" "$LOGSTASH_CONF"
Restart=on-failure
WorkingDirectory=/
PrivateTmp=true

[Install]
WantedBy=multi-user.target
