# lokinet-default-config.service
[Unit]
Description=Lokinet default config generator
ConditionPathExists=!/var/lib/lokinet/lokinet.ini
Wants=local-fs.target
After=local-fs.target

[Service]
Type=oneshot
User=lokinet
WorkingDirectory=/var/lib/lokinet
# at the time of writing it generates wrong paths when invoked without explicit path
ExecStart=/usr/bin/lokinet -g ./lokinet.ini
