# loolwsd.service
[Unit]
Description=loolwsd as a service

[Service]
User=lool
Group=lool
ExecStart=/usr/bin/loolwsd 
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
