# lsi_mrdsnmpd.service
[Unit]
Description=LSI MegaRAID Storage Manager SMNP Agent
Requires=snmpd.service lsi_msm.service

[Service]
Type=forking
EnvironmentFile=/etc/profile.d/msm.sh
ExecStart=/usr/bin/lsi_mrdsnmpagent -c /etc/snmp/snmpd.conf

[Install]
WantedBy=multi-user.target
