# lsyncd.service
[Unit]
Description=synchronization daemon

[Service]
Type=simple
User=lsyncd
Group=lsyncd
Restart=always
ExecStart=/usr/bin/lsyncd -nodaemon /etc/lsyncd/lsyncd.conf.lua
ExecReload=/bin/kill -HUP $MAINPID

# Hardening
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
