# luky-borg-backup.timer
[Unit]
Description=Run luky-borg-backup every night

[Timer]
OnCalendar=04:00
AccuracySec=1h
Persistent=yes

[Install]
WantedBy=timers.target
