# luky-borg-backup.service
[Unit]
Description=LukyLX Incremental Backup

[Service]
Type=oneshot
ExecStart=/usr/bin/luky-borg-backup
