# luky-borg-backup.service
[Unit]
Description=LukyLX Incremental Backup
Wants=unit-status-mail@%n.service
Before=unit-status-mail@%n.service

[Service]
Type=oneshot
ExecStart=/usr/bin/luky-borg-backup
