# lxd.socket
[Unit]
Description=LXD - unix socket
Documentation=file:///usr/share/doc/lxd/

[Socket]
ListenStream=/var/lib/lxd/unix.socket
SocketMode=0660
SocketGroup=lxd
Service=lxd.service

[Install]
WantedBy=sockets.target

