# macfanctld.service
[Unit]
Description=Mac Fan Control Daemon
Documentation=man:macfanctld(1)

[Service]
ExecStart=/usr/bin/macfanctld -f

[Install]
WantedBy=multi-user.target
