# maestral@.service
[Unit]
Description = Maestral daemon for config %i
Documentation = https://maestral.app

[Service]
Type = notify
NotifyAccess = exec
ExecStart = /usr/bin/maestral start -f -c %i
ExecStop = /usr/bin/maestral stop -c %i
ExecStopPost=/usr/bin/env bash -c "if [ ${SERVICE_RESULT} != success ];\
then notify-send Maestral 'Daemon failed'; fi"
WatchdogSec = 30s

[Install]
WantedBy = default.target
