# maestral.service
[Unit]
Description = Maestral daemon
Documentation = "https://github.com/SamSchott/maestral-dropbox/wiki"

[Service]
Type = notify
NotifyAccess = exec
ExecStart = /usr/bin/maestral start --foreground
ExecReload = /usr/bin/maestral restart
ExecStop = /usr/bin/maestral stop
RestartSec = 1
Restart = on-failure
WatchdogSec = 30s

[Install]
WantedBy = default.target
