# mage-server.service
[Unit]
Description=Starts mage-server on boot 

[Service]
Type=simple
ExecStart=/usr/bin/mage-server

[Install]
WantedBy=multi-user.target
