# mailcatcher-smtp.socket
[Unit]
Description=Mailcatcher SMTP Activation Socket
After=network.target

[Socket]
ListenStream=25

[Install]
WantedBy=sockets.target
