# mailcatcher-http.socket
[Unit]
Description=Mailcatcher HTTP Activation Socket
Requires=mailcatcher.service
After=network.target mailcatcher.service

[Socket]
ListenStream=80
Accept=true

[Install]
WantedBy=sockets.target
