# mailcatcher-http.socket
[Unit]
Description=Mailcatcher HTTP Activation Socket
After=network.target

[Socket]
ListenStream=80

[Install]
WantedBy=sockets.target
