# mailcatcher-smtp.socket
[Unit]
Description=Mailcatcher SMTP Activation Socket
Requires=mailcatcher.service
After=network.target mailcatcher.service

[Socket]
ListenStream=25
Accept=true

[Install]
WantedBy=sockets.target
