# mailcatcher-http.socket
[Unit]
Description=Mailcatcher HTTP Activation Socket

[Socket]
ListenStream=80
Accept=true

[Install]
WantedBy=sockets.target
