# mailpile.service
[Unit]
Description=Mailpile Daemon
Requires=network.target
After=network.target

[Service]
User=mailpile
PIDFile=/run/mailpile.pid
PrivateTmp=true
Environment=MAILPILE_HOME=/var/lib/mailpile
ExecStart=/usr/bin/mp --www= --wait

[Install]
WantedBy=multi-user.target
