# mailcatcher.service
[Unit]
Description=Mailcatcher Service
After=network.target

[Service]
ExecStart=/usr/bin/mailcatcher --ip localhost --http-port 1080 --smtp-port 1025 --foreground --verbose
ExecStartPost=/usr/bin/sleep 1
ExecStop=/usr/bin/curl --silent --get --request DELETE localhost:1080

[Install]
WantedBy=multi-user.target
