# mailcatcher-http@.service
[Unit]
Description=Mailcatcher HTTP Mapping Service

[Service]
ExecStart=/usr/bin/nc localhost 1080
StandardInput=socket
StandardOutput=socket
