# mailcatcher-smtp.service
[Unit]
Description=Mailcatcher SMTP Mapping Service
Requires=mailcatcher.service
After=mailcatcher.service

[Service]
ExecStart=/usr/lib/systemd/systemd-socket-proxyd localhost:1025
