# maldet-update-signatures.timer
[Unit]
Description=Linux Malware Detect signature update timer - maldet

[Timer]
OnCalendar=daily
RandomizedDelaySec=360
Persistent=true

[Install]
WantedBy=timers.target
