# mandoc.timer
[Unit]
Description=Daily manpage database update

[Timer]
OnCalendar=daily
AccuracySec=1d
Persistent=true
