# manolo.luksopenwithoutkey.service
# selbst erstellt 20.1.18

[Unit]
Description=Unlock Encrypted Device for MD RAID 1
Conflicts=umount.target
DefaultDependencies=no
BindsTo=dev-mapper-%i.device
Before=mnt-daten.mount
After=dev-md1.device mdmonitor.service
Requires=dev-md1.device mdmonitor.service
ConditionPathExists=!/dev/disk/by-uuid/BA19-64EE

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/lib/systemd/systemd-cryptsetup attach 'storage' '/dev/md1' 'none' 'luks';

[Install]
WantedBy=basic.target
