# maradns-zoneserver.service
[Unit]
Description=MaraDNS Zoneserver handles DNS zone transfers and any TCP DNS queries
Documentation=man:zoneserver(8)
Requires=network.target maradns.service
After=network.target

[Service]
ExecStart=/usr/sbin/zoneserver -f /etc/maradns/mararc
StandardOutput=syslog
Restart=always

[Install]
WantedBy=multi-user.target
