# markdump.service
[Unit]
Description=markdump
After=network.target

[Service]
Type=simple
User=markdump
Group=markdump
PrivateDevices=true
PrivateIPC=true
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectSystem=strict
ReadWritePaths=/var/lib/markdump
RuntimeDirectory=markdump
EnvironmentFile=/etc/markdump/markdump.conf
ExecStart=/usr/bin/markdump

[Install]
WantedBy=multi-user.target
