# mastodon-sidekiq.service
[Unit]
Description=Mastodon background queue service
After=network.target

[Service]
Type=simple
User=mastodon
WorkingDirectory=/var/lib/mastodon
Environment="RAILS_ENV=production"
Environment="DB_POOL=5"
ExecStart=/usr/bin/bundle exec sidekiq -c 5 -q default -q mailers -q pull -q push
TimeoutSec=15
Restart=always

[Install]
WantedBy=multi-user.target
