# mastodon-web.service
[Unit]
Description=Mastodon web workers service
After=network.target

[Service]
Type=simple
User=mastodon
WorkingDirectory=/var/lib/mastodon
Environment="RAILS_ENV=production"
Environment="PORT=3000"
ExecStart=/usr/bin/bundle exec puma -C config/puma.rb
; ExecReload=/bin/kill -SIGUSR1 $MAINPID
TimeoutSec=15
Restart=always

[Install]
WantedBy=multi-user.target
