# matrix-puppet-hangouts@.service
# /usr/lib/systemd/system/matrix-puppet-hangouts@.service
# /etc/systemd/system/hangouts-puppet.service
[Unit]
Description=Hangouts Matrix Puppet

[Service]
Type=simple
Environment=DEBUG=verbose:matrix-puppet:*
Environment=HOME=/var/lib/matrix-puppet-hangouts/%i/
User=matrix-hangouts
Group=matrix-hangouts
WorkingDirectory=/var/lib/matrix-puppet-hangouts/%i/
ExecStart=/usr/bin/node /usr/lib/node_modules/matrix-puppet-hangouts/index.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
