# matrix-registration.service
[Unit]
Description=A token based matrix registration api.
After=network-online.target
Requires=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/matrix_registration --config-path /etc/matrix-registration/config.yaml serve
WorkingDirectory=/var/lib/matrix-registration
DynamicUser=true
ConfigurationDirectory=matrix-registration
StateDirectory=matrix-registration

[Install]
WantedBy=multi-user.target
