# matterircd.service
[Unit]
Description = minimal IRC server which integrates with mattermost and slack

[Service]
ExecStart = /usr/bin/matterircd --conf /etc/matterircd/config.toml
User = matterircd

[Install]
WantedBy = multiuser.target
