# maubot.service
[Unit]
Description=Plugin-based Matrix bot

[Service]
Type=simple
User=maubot
Group=maubot
WorkingDirectory=/var/lib/maubot/
ExecStart=/bin/python -m maubot -c /etc/maubot/config.yaml
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
