# maubot.service
[Unit]
Description=Plugin-based Matrix bot

[Service]
Type=simple
User=maubot
Group=maubot
WorkingDirectory=/var/lib/maubot/
ExecStart=/usr/bin/maubot -c /etc/synapse/maubot/config.yaml \
          -r /etc/synapse/maubot/registration.yaml
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
