# mautrix-telegram.service
[Unit]
Description=A Matrix-Telegram hybrid puppeting/relaybot bridge

[Service]
Environment=FFMPEG_BINARY=/usr/bin/ffmpeg
ExecStart=/usr/bin/mautrix-telegram -b /etc/mautrix-telegram/config.yaml -c /etc/mautrix-telegram/config.yaml -r /etc/mautrix-telegram/registration.yaml
WorkingDirectory=/etc/mautrix-telegram/
Restart=on-failure

[Install]
WantedBy=multi-user.target
