# mautrix-wsproxy.service
[Unit]
Description=Mautrix wsproxy

[Service]
ExecStart=/usr/bin/mautrix-wsproxy -config /etc/mautrix-wsproxy/wsproxy.yaml
Restart=on-failure
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes

[Install]
WantedBy=multi-user.target
