# mautrix-googlechat.service
[Unit]
Description=Mautrix Google Chat bridge

[Service]
ExecStart=/usr/bin/python -m mautrix_googlechat -c /etc/mautrix-googlechat/config.yaml

[Install]
WantedBy=multi-user.target
