# mcproxy.service
[Unit]
Description=mcproxy multicast proxy for IGMP/MLD (RFC 4605)
Documentation=http://mcproxy.realmv6.org/wiki/Documentation
After=network.target pihole-FTL.service iptables.service

[Service]
ExecStart=/usr/bin/mcproxy -r -f /etc/mcproxy/mcproxy.conf

[Install]
WantedBy=multi-user.target
