# mdm-plymouth.service
[Unit]
Description=LightDM Display Manager
Documentation=man:mdm(1)
After=systemd-user-sessions.service getty@tty1

Conflicts=plymouth-quit.service
After=plymouth-quit.service

[Service]
ExecStart=/usr/sbin/mdm --nodaemon
ExecReload=/usr/sbin/mdm-safe-restart
ExecStop=/usr/sbin/mdm-stop
Restart=always
RestartSec=1s
TimeoutStopSec=5s
IgnoreSIGPIPE=no

[Install]
WantedBy=graphical.target
Alias=display-manager.service
