# me176c.mount
[Unit]
Description=ME176C Data Partition

[Mount]
What=/dev/disk/by-partlabel/config
Where=/me176c
Type=ext4
Options=nosuid,nodev,noatime,data=ordered
