# mediatomb.service
[Unit]
Description=MediaTomb Daemon
After=network.target

[Service]
EnvironmentFile=/etc/conf.d/mediatomb
ExecStart=/usr/bin/mediatomb --user $MT_USER --group $MT_GROUP --port $MT_PORT --home $MT_HOME --logfile $MT_LOGFILE --pidfile $MT_PIDFILE $MT_OPTIONS
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target