# mercure.service
[Unit]
Description=Mercure - Server-sent live updates
After=network.target

[Service]
ExecStart=/usr/bin/mercure start
User=mercure
Group=mercure
RuntimeDirectory=mercure

[Install]
WantedBy=multi-user.target
