# midnight-dler.service
[Unit]
Description=Update and Shutdown Service
After=network-online.target
#SuccessAction=none
#OnSuccess=systemd-suspend.service
[Service]
RemainAfterExit=no
Type=oneshot
ExecStart=/usr/bin/midnight-dler.sh 
[Install]
WantedBy=default.target
