# milter-greylist.service
[Unit]
Description=Greylist milter
Before=postfix.service
After=network.target

[Service]
UMask=0007
Environment=SOCKET=/var/lib/milter-greylist/greylist.sock
ExecStart=/usr/bin/milter-greylist -D -d /var/lib/milter-greylist/greylist.db -f /etc/mail/greylist.conf -u milt-gl:milt-gl -p $SOCKET

[Install]
WantedBy=multi-user.target
