# milter-spamd.service
[Unit]
Description = Milter for SpamAssassin
Wants = spamassassin.service
After=network.target spamassassin.service

[Service]
Type = forking
User = milter-spamd
Group = mail
RuntimeDirectory=milter-spamd
RuntimeDirectoryMode=0710
ExecStart = /usr/bin/milter-spamd -p /run/milter-spamd/milter-spamd.sock -U spamd

[Install]
WantedBy = multi-user.target
