# milter-greylist.service
[Unit]
Description=Greylist milter
Before=postfix.service
After=network.target

[Service]
UMask=0007
ExecStart=/usr/bin/milter-greylist -D

[Install]
WantedBy=multi-user.target
