# minecraftd.service
[Unit]
Description=Minecraft Server
After=local-fs.target network.target

[Service]
Type=forking
ExecStart=/usr/bin/minecraftd start
ExecStop=/usr/bin/minecraftd stop
USER=minecraft
GROUP=minecraft

[Install]
WantedBy=multi-user.target
