# minecraftd.service
[Unit]
Description=Minecraft Server
After=local-fs.target network.target

[Service]
Type=forking
EnvironmentFile=/etc/conf.d/minecraft
ExecStart=/usr/bin/minecraftd start
ExecStop=/usr/bin/minecraftd stop
USER=$MC_USER
GROUP=$MC_USER

[Install]
WantedBy=multi-user.target
