# minisatip.service
[Unit]
Description=Mini SAT>IP server
After=network-online.target
Wants=network-online.target oscam.service

[Service]
User=minisatip
Group=minisatip
Type=simple
AmbientCapabilities=CAP_NET_BIND_SERVICE
RuntimeDirectory=libmcli
EnvironmentFile=-/etc/conf.d/minisatip
ExecStart=/usr/bin/minisatip -f -R /usr/share/minisatip/html $SERVER_ARGS

[Install]
WantedBy=multi-user.target
