# miniflux.service
[Unit]
Description=Miniflux Feed Reader
After=network.target
After=postgresql.service

[Service]
Type=simple
EnvironmentFile=/etc/miniflux.conf
User=miniflux
ExecStart=/usr/bin/miniflux
Restart=always

[Install]
WantedBy=multi-user.target
