# minimserver@.service
[Unit]
Description=MinimServer

[Service]
User=%I
RemainAfterExit=true
ExecStart=/opt/minimserver/bin/startd
ExecStop=/opt/minimserver/bin/stopall

[Install]
WantedBy=network.target