# miredo.service
[Unit]
Description=Teredo IPv6 tunneling
After=network.target

[Service]
ExecStartPre=/usr/bin/miredo-checkconf -f /etc/miredo/miredo.conf
ExecStart=/usr/bin/miredo -f
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
