# mlink.service
[Unit]
Description=Start mlink bridge to connect mBlock web version to device

[Service]
Type=simple
ExecStart=/usr/bin/mblock-mlink start
Restart=always
RestartSec=30

[Install]
WantedBy=default.target
