# mmonit.service
[Unit]
Description = Easy, proactive monitoring of Unix systems, network and cloud services
After = network.target

[Service]
Type=forking
ExecStart = /opt/mmonit/bin/mmonit start
ExecStop = /opt/mmonit/bin/mmonit stop
PIDFile = /opt/mmonit/logs/mmonit.pid

[Install]
WantedBy = multi-user.target
