# modacom-u3500.service
[Unit]
Description=Mezon modem daemon
After=network.target

[Service]
ExecStart=/usr/bin/gctwimax -C /usr/share/gctwimax/gctwimax.conf
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target
