# molly-brown.service
[Unit]
Description=Molly Brown gemini server
After=network.target

[Service]
Type=simple
NoNewPrivileges=true
ProtectSystem=true
ProtectHome=true
PrivateDevices=true
User=molly
ExecStart=/usr/bin/molly-brown -c /etc/molly.conf

[Install]
WantedBy=multi-user.target
