# monica-scheduler.timer
[Unit]
Description=Monica scheduler timer

[Timer]
# every ten minutes
OnCalendar=*:0/10

[Install]
WantedBy=timers.target
