# moolticute.service
[Unit]
Description=Moolticute daemon
Type=Simple

[Service]
ExecStart=/usr/bin/moolticuted
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target
