# moolticute.service
[Unit]
Description=Moolticute daemon

[Service]
ExecStart=/usr/bin/moolticuted
KillMode=process
Restart=always
Type=simple

[Install]
WantedBy=multi-user.target
