# moonraker.service
[Unit]
Description=Moonraker Klipper HTTP Server
Requires=klipper.service
After=network.target klipper.service

[Service]
Type=simple
User=klipper
SupplementaryGroups=moonraker-admin
SyslogIdentifier=moonraker
RemainAfterExit=yes
EnvironmentFile=/var/opt/moonraker/systemd/moonraker.env
ExecStart=/usr/bin/python /opt/moonraker/moonraker/moonraker.py -n $MOONRAKER_ARGS
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
