# moonshine.service
[Unit]
Description=Streaming server using the NVIDIA GameStream / Moonlight protocol.
PartOf=graphical-session.target
Wants=graphical-session.target
After=graphical-session.target

[Service]
ExecStart=/usr/bin/moonshine /etc/moonshine/config.toml -v
WorkingDirectory=/etc/moonshine

[Install]
WantedBy=graphical-session.target
