# moonshine.service
[Unit]
Description=Streaming server using the NVIDIA GameStream / Moonlight protocol.
PartOf=graphical-session.target
Wants=graphical-session.target
After=graphical-session.target

[Service]
ExecStart=/usr/bin/moonshine %E/moonshine/config.toml

[Install]
WantedBy=graphical-session.target
