# mopidy@.service
[Unit]
Description=mopidy
After=network.target

[Service]
User=%I
ExecStart=/usr/bin/mopidy

[Install]
WantedBy=multi-user.target
