# mopidy-user.service
[Unit]
Description=mopidy
After=network.target

[Service]
ExecStart=/usr/bin/mopidy

[Install]
WantedBy=default.target
