# mosix.service
[Unit]
Description=Mosix Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/sbin/mosd

[Install]
WantedBy=multi-user.target
