# motioneye.service
[Unit]
Description=motionEye server
After=local-fs.target network.target

[Service]
User=motion
Environment="LANGUAGE=en_EN"
RuntimeDirectory=motioneye
LogsDirectory=motioneye
StateDirectory=motioneye
ExecStart=/usr/bin/meyectl startserver -c /etc/motioneye/motioneye.conf
Type=simple
Restart=on-abort

[Install]
WantedBy=multi-user.target
