# movim.service
[Unit]
Description=Movim Daemon
After=network.target

[Service]
User=movim
WorkingDirectory=/usr/share/webapps/movim
EnvironmentFile=/etc/default/movim
ExecStartPre=/usr/bin/php vendor/bin/phinx migrate
ExecStart=/usr/bin/php daemon.php start --url=${MOVIM_URL} --port=${MOVIM_PORT} --interface=${MOVIM_INTERFACE}
Restart=always

[Install]
WantedBy=multi-user.target
