# mpd-user.socket
[Unit]
Description=Music Player Daemon socket

[Socket]
ListenStream=6600
Service=mpd.service

[Install]
WantedBy=default.target
