# mpd-system.socket
[Unit]
Description=Music Player Daemon socket
After=sound.target

[Socket]
ListenStream=6600
Service=mpd.service

[Install]
WantedBy=sockets.target
