# mpdrandom-user.service
[Unit]
Description=Play a random album in mpd
After=mpd.service

[Service]
ExecStart=/usr/bin/mpdrandom -d

[Install]
WantedBy=default.target
