# mrouted.service
[Unit]
Description=Multicast Route Daemon
After=network.target
ConditionPathExists=/etc/mrouted.conf
Conflicts=pimd.service smcroute.service

[Service]
Type=forking
PIDFile=/var/run/mrouted.pid
ExecStart=/usr/bin/mrouted
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-abort

[Install]
WantedBy=multi-user.target
