# msmtp.service
[Service]
Type=oneshot
ExecStart=/usr/bin/msmtp-flush-mail %h/.local/share/msmtp-queue
Restart=on-failure
RestartSec=1min
