# msmtp.path
[Path]
DirectoryNotEmpty=%h/.local/share/msmtp-queue
MakeDirectory=yes
DirectoryMode=0700

[Install]
WantedBy=paths.target
