# mstpd.service
[Unit]
Description=Multiple Spanning Tree daemon

[Service]
Type=forking
ExecStart=/usr/bin/mstpd

[Install]
WantedBy=multi-user.target
