# mstream.service
[Unit]
Description=mstream
After=network.target

[Service]
User=mstream
Type=simple
WorkingDirectory=/var/lib/mstream
ExecStart=/usr/bin/mstream -j /var/lib/mstream/config.json

[Install]
WantedBy=multi-user.target

