# mstream.service
[Unit]
Description=mStream
After=network.target

[Service]
User=mstream
Type=simple
WorkingDirectory=/var/lib/mstream
ExecStart=/usr/bin/mstream -j /etc/mstream.json

[Install]
WantedBy=multi-user.target
