# mstream.service
[Unit]
Description=mstream
After=network.target

[Service]
User=mstream
Type=simple
WorkingDirectory=/var/lib/mstream
ExecStart=/usr/bin/mstream /var/lib/mstream/config.json -m /var/lib/mstream/music

[Install]
WantedBy=multi-user.target

