# mstream.service
[Unit]
Description=mStream Daemon
After=network.target

[Service]
Type=simple

ExecStart=/usr/bin/mstream --json=/etc/mstream.json
WorkingDirectory=/var/lib/mstream
Restart=on-success

User=mstream
Group=mstream

[Install]
WantedBy=multi-user.target
