# mtg.service
[Unit]
Description=Bullshit-free MTProto proxy for Telegram
After=network-online.target

[Service]
User=nobody
Group=nobody
EnvironmentFile=/etc/conf.d/mtg
ExecStart=/usr/bin/mtg $SECRET $ADTAG

[Install]
WantedBy=multi-user.target

