# mtg@.service
[Unit]
Description=Bullshit-free MTProto proxy for Telegram
After=network-online.target

[Service]
DynamicUser=yes
ExecStart=/usr/bin/mtg run /etc/mtg/%i.toml

[Install]
WantedBy=multi-user.target
