# mtproxy-config.timer
[Unit]
Description=Daily update of Telegram configuration

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
