# mullvadd.initd
#!/sbin/openrc-run
description="Mullvad VPN Service"

depend() {
    need net
}

supervisor="supervise-daemon"
command="/opt/Mullvad\ VPN/resources/mullvad-daemon"
command_args="${MULLVPN_OPTS}"
#pidfile="/run/${SVCNAME}.pid"
command_background=true
