# murmur.service
[Unit]
Description=Mumble Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/mumble-server -ini /etc/murmur.ini -fg
Restart=always

[Install]
WantedBy=multi-user.target
