# mutekey-fix.service
[Unit]
Descrition=GPD Win Max mute key fix

[Service]
ExecStart=/usr/local/bin/mutekey-fix

[Install]
WantedBy=multi-user.target
