# mycroft.service
[Unit]
Description=Mycroft main service
After=network.target local-fs.target remote-fs.target

[Service]
User=mycroft
Group=mycroft
SyslogIdentifier=mycroft-core
Type=forking
RemainAfterExit=no
ExecStart=/usr/share/mycroft-core/start-mycroft.sh all
ExecStop=/usr/share/mycroft-core/stop-mycroft.sh all
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
