# mycroft-voice.service
[Unit]
Description=Mycroft voice processing
After=network.target local-fs.target remote-fs.target mycroft-skills

[Service]
Type=simple
PermissionsStartOnly=true
ExecStart=/bin/bash -c "source /usr/share/mycroft-core/.virtualenvs/mycroft/bin/activate && /usr/share/mycroft-core/start.sh voice"

[Install]
WantedBy=multi-user.target
