# mycroft-service.service
[Unit]
Description=Mycroft main service
After=network.target local-fs.target remote-fs.target

[Service]
Type=simple
PermissionsStartOnly=true
ExecStart=/bin/bash -c "source /usr/share/mycroft-core/.virtualenvs/mycroft/bin/activate && /usr/share/mycroft-core/start.sh service"

[Install]
WantedBy=multi-user.target
