# mycron@.service
[Unit]
Description=mycron for %i
Wants=multi-user.target

[Service]
User=%i
Restart=always
ExecStart=/usr/bin/mycron

[Install]
WantedBy=multi-user.target
