# mylar.service
[Unit]
Description=Comic Book add-on for SABnzbd+

[Service]
ExecStart=/usr/bin/python2 /usr/share/mylar/Mylar.py --nolaunch --config /etc/mylar/config.ini --datadir /var/lib/mylar

[Install]
WantedBy=multi-user.target
