# mysql-router.service
[Unit]
Description=MySQL router
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/mysqlrouter --config /etc/mysql/router.ini
ProtectSystem=full
ProtectHome=on
PrivateDevices=on
NoNewPrivileges=on

[Install]
WantedBy=multi-user.target
