# mythbackend.service
[Unit]
Description=MythTV Backend
Wants=network-online.target mysqld.service systemd-udev-settle.service
After=network.target mysqld.service systemd-udev-settle.service

[Service]
Type=simple
User=mythtv
ExecStart=/usr/bin/mythbackend
Restart=on-failure

[Install]
WantedBy=multi-user.target
