# nagios.service
[Unit]
Description=Nagios Monitoring System
After=network.target

[Service]
ExecStart=/usr/bin/nagios /etc/nagios/nagios.cfg
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process

[Install]
WantedBy=multi-user.target
