# namecoin.service
[Unit]
Description=Namecoin Daemon
After=network.target

[Service]
User=namecoin
Group=namecoin
ExecStart=/usr/bin/namecoind -daemon=0 \
                             -conf=/etc/namecoin/namecoin.conf \
                             -datadir=/var/lib/namecoin \
                             -pid=/run/namecoin/namecoin.pid
ExecReload=/usr/bin/kill -HUP $MAINPID
ExecStop=/usr/bin/namecoind stop

[Install]
WantedBy=multi-user.target
