# namecoind@.service
[Unit]
Description=Namecoin Core Daemon
After=local-fs.target network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/namecoind
ExecStop=/usr/bin/namecoin-cli stop
User=%I

[Install]
WantedBy=default.target

