# namecoin.service
[Unit]
Description=Namecoin Daemon
After=network.target

[Service]
User=namecoin
Group=namecoin
ExecStart=/usr/bin/namecoind -daemon=0                       \
                            -conf=/etc/namecoin/namecoin.conf \
                            -datadir=/srv/namecoin           \
                            -pid=/run/namecoind.pid
ExecReload=/usr/bin/kill -HUP $MAINPID
ExecStop=/usr/bin/namecoind stop

[Install]
WantedBy=multi-user.target
