# nbxplorer.service
[Unit]
Description=NBXplorer
After=network.target
After=bitcoind.service

[Service]
Type=oneshot
ExecStart=/usr/bin/nbxplorer
RemainAfterExit=yes
Restart=on-failure
#User=

[Install]
WantedBy=multi-user.target
