# ncdns.service
[Unit]
Description=Namecoin to DNS bridge daemon

[Service]
# If you don't want to bind to port 53, then uncomment to run as non-root:
# User=ncdns
# Group=ncdns
ExecStart=/usr/bin/ncdns
Restart=on-failure

[Install]
WantedBy=default.target
