# ndyndns.service
[Unit]
Description=Dynamic DNS updater
After=network.target

[Service]
ExecStart=/usr/bin/ndyndns -n -f /etc/ndyndns.conf

[Install]
WantedBy=network.target
