# nemo-desktop.service
[Unit]
Description=Nemo desktop
DefaultDependencies=false
After=dbus.socket
Before=basic.target shutdown.target
Conflicts=shutdown.target

[Service]
ExecStart=/usr/bin/nemo-desktop
