# neo4j.service
[Unit]
Description=Neo4j

[Service]
User=neo4j
Type=forking
ExecStart=/usr/bin/neo4j start
ExecStop=/usr/bin/neo4j stop
PIDFile=/run/neo4j/neo4j-service.pid
