# neofetch-port.service
[Unit]
Description=Run neofetch-port

[Service]
ExecStart=/usr/bin/neofetch-port
Restart=on-failure

[Install]
WantedBy=network.target
