# netbox-system.service
[Unit]
Description=Netbox Daemon
After=network.target

[Service]
User=netbox
Type=simple
PIDFile=/run/netbox/netbox.pid
ExecStart=/usr/bin/gunicorn -c /etc/netbox/gunicorn_config.py --pythonpath /opt/netbox -p /run/netbox/netbox.pid netbox.wsgi

[Install]
WantedBy=multi-user.target
