# netboard.service
[Unit]
Description=netboard

[Service]
ExecStart=/usr/bin/netboard listen
RestartSec=1s
Restart=on-failure

[Install]
WantedBy=default.target
