# netbox-system.service
[Unit]
Description=Netbox Daemon
After=network.target

[Service]
User=netbox
Type=forking
PIDFile=/run/netbox/netbox.pid
ExecStart=/usr/bin/gunicorn -c /etc/netbox/gunicorn_config.py netbox.wsgi

[Install]
WantedBy=multi-user.target
