# nethogs.service
[Unit]
Description=nethogs tracemode service
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/usr/bin/nethogst

# Hardening
ProtectSystem=strict
ReadWritePaths=/var/lib
StateDirectory=nethogs
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectKernelModules=yes
PrivateTmp=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
RestrictNamespaces=yes

[Install]
WantedBy=multi-user.target
