# netkit-telnetd-ssl@.service
[Unit]
Description=Netkit Telnet Server Instance

[Service]
EnvironmentFile=-/etc/default/telnetd
ExecStart=-/usr/bin/in.telnetd-ssl -z ssl $TELNETD_ARGS
StandardInput=socket
