# networktablet.service
[Unit]
Description=Multimedia Service

[Service]
Type=simple
ExecStart=/usr/bin/networktablet
Restart=on-failure
WantedBy=default.target
