# networkaudio.service
[Unit]
Description=Sygnalist Network Audio Daemon
After=network.target sound.target

[Service]
ExecStart=/usr/bin/networkaudiod

[Install]
WantedBy=multi-user.target
