# networkaudio@.service
[Unit]
Description=Sygnalist Network Audio Daemon
After=network.target sound.target

[Service]
User=%I
ExecStart=/usr/bin/networkaudiod

[Install]
WantedBy=multi-user.target
