# networkmanager-vpn-web-ui.service
[Unit]
Description=Web UI for NetworkManager-managed VPN connections
After=NetworkManager.service

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/networkmanager-vpn-web-ui -bind-address :9092

[Install]
WantedBy=multi-user.target
