# neutron-l3-agent.service
[Unit]
Description=Neutron l3 agent

[Service]
ExecStart=/usr/bin/neutron-l3-agent --config-file=/etc/neutron/neutron.conf --config-file=/etc/neutron/l3_agent.ini --log-file=/var/log/neutron/l3-agent.log
User=neutron

[Install]
WantedBy=multi-user.target
