# neutron-lbaasv2-agent.service
[Unit]
Description=Neutron LBaaSv2 Agent

[Service]
ExecStart=/usr/bin/neutron-lbaasv2-agent --config-file=/etc/neutron/neutron.conf --config-file=/etc/neutron/lbaas_agent.ini --log-file=/var/log/neutron/lbaasv2_agent.log
User=neutron

[Install]
WantedBy=multi-user.target
