# neutron-dhcp-agent.service
[Unit]
Description=Neutron DHCP agent

[Service]
ExecStart=/usr/bin/neutron-dhcp-agent --config-file=/etc/neutron/neutron.conf --config-file=/etc/neutron/dhcp_agent.ini --log-file=/var/log/neutron/dhcp-agent.log
User=neutron

[Install]
WantedBy=multi-user.target
