# neutron-plugin-oneconvergence-agent.service
[Unit]
Description=Neutron One Convergence Plugin agent

[Service]
ExecStart=/usr/bin/neutron-nvsd-agent --config-file=/etc/neutron/neutron.conf --config-file=/etc/neutron/plugins/oneconvergence/nvsdplugin.ini --log-file=/var/log/neutron/nvsd-agent.log
User=neutron

[Install]
WantedBy=multi-user.target
