# neutron-plugin-mlnx-agent.service
[Unit]
Description=Neutron MLNX Plugin Agent

[Service]
ExecStart=/usr/bin/neutron-mlnx-agent --config-file=/etc/neutron/neutron.conf --config-file=/etc/neutron/plugins/mlnx_conf.ini --log-file=/var/log/neutron/mlnx_conf.log
User=neutron

[Install]
WantedBy=multi-user.target
