# neutron-plugin-ryu-agent.service
[Unit]
Description=Neutron ryu plugin agent

[Service]
ExecStart=/usr/bin/neutron-ryu-agent --config-file=/etc/neutron/neutron.conf --config-file=/etc/neutron/plugins/ryu/ryu.ini --log-file=/var/log/neutron/ryu-agent.log
User=neutron

[Install]
WantedBy=multi-user.target
