# neutron-plugin-ibm-agent.service
[Unit]
Description=Neutron IBM SDN agent

[Service]
ExecStart=/usr/bin/neutron-ibm-agent --config-file=/etc/neutron/neutron.conf --config-file=/etc/neutron/plugins/plugins/ibm/sdnve_neutron_plugin.ini --log-file=/var/log/neutron/ibm-agent.log
User=neutron

[Install]
WantedBy=multi-user.target
