# nextcloud-news-updater.service
[Unit]
After=default.target

[Service]
Type=simple
User=http
ExecStart=/usr/bin/owncloud-news-updater -c /etc/webapps/nextcloud/news/nextcloud-news-updater.ini

[Install]
WantedBy=default.target
