# nextcloud-preview-generate-all.timer
[Unit]
Description=Nextcloud; weekly generate previews for all files

[Timer]
OnBootSec=20min
OnCalendar=weekly
Persistent=True

[Install]
WantedBy=timers.target