# nextcloud-preview.timer
[Unit]
Description=Nextcloud preview generation

[Timer]
OnBootSec=10min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target