# nextcloud-app-facerecognition.timer
[Unit]
Description=Run Nextcloud FaceRecognition every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
