# nextcloud-files-scan-all.timer
[Unit]
Description=Nextcloud; weekly scan all files

[Timer]
OnBootSec=15min
OnCalendar=weekly
Persistent=True

[Install]
WantedBy=timers.target