# nextcloud-cron.timer
[Unit]
Description=Nextcloud; run cron.php periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target