# nextcloud-cron.timer
[Unit]
Description=Nextcloud; periodically run cron.php

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target