# nexus.service
[Unit]
Description=Sonatype Nexus 2.x
After=network.target

[Service]
Environment=PLEXUS_NEXUS_WORK=/var/lib/nexus
ExecStart=/usr/lib/nexus/bin/nexus start
ExecStop=/usr/lib/nexus/bin/nexus stop
ExecReload=/usr/lib/nexus/bin/nexus restart
PIDFile=/run/nexus/nexus.pid
Type=forking
User=nexus
Group=nexus

[Install]
WantedBy=multi-user.target

