# nft-blacklist.service
[Unit]
Description=Blacklist country specific IP blocks using nftables
Requires=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/conf.d/nft-blacklist
ExecStart=/usr/bin/nft-blacklist $NFT_BLACKLIST_ARGS
