# nginx-totp-auth.service
[Unit]
Description=TOTP auth daemon for nginx

[Service]
RuntimeDir=totp_auth
ExecStart=/usr/bin/spawn-fcgi -u http -s /run/totp_auth/sock -M 666 -n /usr/lib/nginx-totp-auth/nginx-totp-auth.bin /etc/nginx-totp-auth.conf

[Install]
WantedBy=multi-user.target
