# nginx.socket
[Unit]
Description=nginx socket

[Socket]
ListenStream=80

[Install]
WantedBy=multi-user.target
