# nginx.service
[Unit]
Description=A high performance web server and a reverse proxy server
After=network.target network-online.target nss-lookup.target

[Service]
Type=forking
PIDFile=nginx.pid
SyslogLevel=err

ExecStart=/usr/lib/systemd/scripts/nginx-user-service
ExecReload=nginx -s reload
KillMode=mixed

[Install]
WantedBy=multi-user.target
