# ninfod.service
[Unit]
Description=IPv6 Node Information Responder

[Service]
Type=forking
ExecStart=/usr/lib/ninfod/ninfod -v -u daemon

[Install]
WantedBy=multi-user.target
