# nis-cachemgr.service
[Unit]
Description=NIS+ cache daemon

[Service]
Type=forking
ExecStart=/usr/bin/nis_cachemgr

[Install]
WantedBy=multi-user.target
