# nitter.service
[Unit]
Description=Nitter Daemon
Requires=redis.service
After=redis.service

[Service]
ExecStart=/usr/bin/nitter

[Install]
WantedBy=multi-user.target
