# nitter.service
[Unit]
Description=Nitter Daemon
After=network.service redis.service
Requires=redis.service

[Service]
User=nitter
ExecStart=/usr/bin/nitter

[Install]
WantedBy=multi-user.target
