# nntp-proxy.service
[Unit]
Description=Simple NNTP proxy
After=network.target

[Service]
User=nntpproxy
ExecStart=/usr/bin/nntp-proxy /etc/nntp-proxy.conf

[Install]
WantedBy=multi-user.target
