# nomad.service
[Unit]
Description=Nomad cluster scheduling agent
Documentation=https://www.nomadproject.io/docs/agent/
Requires=network-online.target
After=network-online.target

[Service]
ExecStart=/usr/bin/nomad agent -config /etc/nomad -data-dir=/var/lib/nomad

[Install]
WantedBy=multi-user.target
