# nomad.service
[Unit]
Description=Nomad agent
Documentation=https://www.nomadproject.io/docs/agent/
Requires=network.target
After=network.target

[Service]
ExecStart=/usr/bin/nomad agent -config /etc/nomad -data-dir=/var/lib/nomad

[Install]
WantedBy=multi-user.target
