# nostromo.service
[Unit]
Description=Nostromo - A simple, fast, and secure HTTP server
After=network.target

[Service]
Type=forking
PIDFile=/usr/share/nostromo/logs/nhttpd.pid
PrivateDevices=yes
ExecStart=/usr/bin/nhttpd
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
