# novacomd.service
[Unit]
Description=Novacom WebOS daemon

[Service]
ExecStart=/usr/bin/novacomd

[Install]
WantedBy=multi-user.target
