# nova-consoleauth.service
[Unit]
Description=Nova ConsoleAuth server

[Service]
ExecStart=/usr/bin/nova-consoleauth --config-dir /etc/nova/
User=nova

[Install]
WantedBy=multi-user.target
