# nova-cert.service
[Unit]
Description=Nova Cert server

[Service]
ExecStart=/usr/bin/nova-cert --config-dir /etc/nova/
User=nova

[Install]
WantedBy=multi-user.target
