# nrclient.service
[Unit]
Description=NeoRouter Client
After=network.target

[Service]
Type=simple
WorkingDirectory=/usr/local/ZebraNetworkSystems/NeoRouter
ExecStart=/usr/bin/nrservice
Restart=always

[Install]
WantedBy=multi-user.target
