# ntopng@.service
[Unit]
Description=ntopng network traffic monitor

[Service]
ExecStart=/usr/local/bin/ntopng -i %I

[Install]
WantedBy=multi-user.target
