# ntopng@.service
[Unit]
Description=ntopng network traffic monitor

[Service]
ExecStart=/usr/bin/ntopng -i eth0

[Install]
WantedBy=multi-user.target
