# nutch-webapp.service
[Unit]
Description=Nutch open source web crawler webapp
After=network.target

[Service]
Type=simple
User=nutch
Group=nutch
PrivateTmp=true
WorkingDirectory=/opt/nutch
ExecStart=/opt/nutch/bin/StartNutchServerWebApp.bash

[Install]
WantedBy=multi-user.target
