# nvidia-fallback.service
[Unit]
Description=Dummy service required by system76-power

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c "echo 'Luke, I am your father' > /dev/null"

[Install]
WantedBy=multi-user.target
