# nvidia-modprobe.service
[Unit]
Description=Unload Nvidia modules before shutdown/reboot.
DefaultDependencies=no
After=final.target

[Service]
Type=oneshot
ExecStart=/usr/bin/modprobe -r nvidia_drm nvidia_modeset nvidia_uvm && /usr/bin/modprobe -r nvidia

[Install]
WantedBy=final.target
