# nvidia-modprobe.service
[Unit]
Description=Ditch nvidia after multi-user.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=/usr/bin/modprobe -r nvidia_drm nvidia_modeset nvidia_uvm nvidia

[Install]
WantedBy=multi-user.target
