# nvidia-modprobe.service
[Unit]
Description=Unload Nvidia modules before shutdown/reboot.

[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=/usr/bin/modprobe -r nvidia_drm nvidia_modeset nvidia_uvm && /usr/bin/modprobe -r nvidia

[Install]
WantedBy=multi-user.target
