# nvx.service
[Unit]
Description="Turn off gpu during boot."

[Service]
Type=oneshot
ExecStart=/usr/bin/nvx off-boot

[Install]
WantedBy=multi-user.target
