# nym-mixnode@.service
[Unit]
Description=Nym Mixnode %I
After=network-online.target nss-lookup.target
Wants=network-online.target nss-lookup.target
Documentation=https://nymtech.net/docs/stable/run-nym-nodes/nodes/mixnodes
StartLimitIntervalSec=350
StartLimitBurst=10

[Service]
User=nym
ExecStart=/usr/bin/nym-mixnode run --id %i
KillSignal=SIGINT
Restart=always
RestartSec=30
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
