# nym-gateway@.service
[Unit]
Description=Nym Gateway %I
After=network-online.target nss-lookup.target
Wants=network-online.target nss-lookup.target
Documentation=https://nymtech.net/docs/stable/run-nym-nodes/nodes/gateways
StartLimitIntervalSec=350
StartLimitBurst=10

[Service]
User=nym
ExecStart=/usr/bin/nym-gateway run --id %i
KillSignal=SIGINT
Restart=on-failure
RestartSec=30
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
