# nzbhydra2.service
[Unit]
Description=NZBHydra2 Daemon
Documentation=https://github.com/theotherp/nzbhydra2
After=network.target

[Service]
User=nzbhydra2
Group=nzbhydra2
Type=simple
Environment="NZBHYDRA_DISABLE_WRAPPER_CHECK=1"
ExecStart=/usr/bin/nzbhydra2 --nobrowser --datafolder /var/lib/nzbhydra2
Restart=always

[Install]
WantedBy=multi-user.target
