# oauth2-proxy@.service
[Unit]
Description=oauth2-proxy server
After=network.target

[Service]
User=oauth2-proxy
Group=oauth2-proxy

Type=simple
ExecStart=/usr/bin/oauth2-proxy --config=/etc/oauth2-proxy/%i.cfg
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict

[Install]
WantedBy=multi-user.target