# octoprint.service
[Unit]
Description=Responsive web interface for controlling a 3D printer
After=network.target

[Service]
User=octoprint
Group=octoprint
Type=simple
ExecStart=/usr/bin/octoprint serve

[Install]
WantedBy=multi-user.target
