# odoo.service
[Unit]
Description=Advanced OpenSource ERP and CRM server
Requires=postgresql.service
After=postgresql.service


[Service]
Type=simple
PermissionsStartOnly=true
EnvironmentFile=-/etc/conf.d/odoo
User=odoo
Group=odoo
SyslogIdentifier=Odoo
PIDFile=/run/odoo/odoo.pid
ExecStartPre=/usr/bin/install -d -m755 -o odoo -g odoo /run/odoo
ExecStart=/usr/bin/odoo.py -c $ODOO_CONF --pidfile=/run/odoo/odoo.pid $ODOO_ARGS
ExecStop=/bin/kill $MAINPID

[Install]
WantedBy=multi-user.target
