# odrive.service
[Unit]
Description=odrive agent
After=local-fs.target network.target

[Service]
Type=simple
ExecStart=/usr/bin/odriveagent
Restart=on-failure

[Install]
WantedBy=multi-user.target
