# offeringd@.service
[Unit]
Description=Offerings Daemon
After=network.target 

[Service]
User=%i
Group=users
Type=simple
ExecStart=/usr/bin/offeringd

[Install]
WantedBy=multi-user.target
