# olaris.service
[Unit]
Description=Olaris Media Server

[Service]
Type=exec
DynamicUser=yes
CacheDirectory=olaris
StateDirectory=olaris
RuntimeDirectory=olaris
Environment=XDG_CACHE_HOME=/var/cache/olaris
Environment=OLARIS_CONFIG_DIR=/var/lib/olaris
WorkingDirectory=/var/lib/olaris
ExecStart=/usr/bin/olaris --verbose=false
Restart=on-failure

[Install]
WantedBy=multi-user.target
