# olivetin.service
[Unit]
Description=OliveTin

[Service]
WorkingDirectory=/usr/share/webapps/olivetin
User=root
ExecStart=/usr/bin/olivetin -configdir /etc/olivetin
Restart=on-failure

[Install]
WantedBy=multi-user.target
