# olivetin@.service
[Unit]
Description=OliveTin run by %I

[Service]
WorkingDirectory=/usr/share/webapps/olivetin
User=%i
ExecStart=/usr/bin/olivetin -configdir /etc/olivetin/%i
Restart=on-failure

[Install]
WantedBy=multi-user.target
