# ombi.service
[Unit]
Description=Ombi
After=network.target

[Service]
User=ombi
Group=ombi
WorkingDirectory=/opt/Ombi
Type=simple
ExecStart=/opt/Ombi/Ombi
TimeoutStopSec=30
Restart=always

[Install]
WantedBy=multi-user.target
