# ombi.service
[Unit]
Description=Ombi
After=network.target

[Service]
User=ombi
Group=ombi
WorkingDirectory=/opt/ombi
Type=simple
ExecStart=/opt/ombi/Ombi
Restart=always
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
