# onedev.service
[Unit]
Description=A self-hosted Git Server with CI/CD and Kanban 
After=syslog.target

[Service]
Type=simple
User=onedev
Group=onedev
ExecStart=/opt/onedev/bin/server.sh console
KillMode=control-group
Environment=SYSTEMD_KILLMODE_WARNING=true

[Install]
WantedBy=multi-user.target
