# onedev.service
[Unit]
Description=A self-hosted Git Server with CI/CD and Kanban 
After=syslog.target

[Service]
Type=forking
User=onedev
Group=onedev
ExecStart=/opt/onedev/bin/server.sh start sysd
ExecStop=/opt/onedev/bin/server.sh stop sysd
KillMode=control-group
Environment=SYSTEMD_KILLMODE_WARNING=true

[Install]
WantedBy=multi-user.target