# onlyoffice-fileconverter.service
[Unit]
Description=Onlyoffice documentserver
After=network.target

[Service]
User=onlyoffice
Type=simple
Environment="NODE_ENV=production-linux" "NODE_CONFIG_DIR=/etc/webapps/onlyoffice/documentserver"
WorkingDirectory=/usr/share/webapps/onlyoffice/documentserver/server/FileConverter/sources/
ExecStart=/usr/bin/node /usr/share/webapps/onlyoffice/documentserver/server/FileConverter/sources/convertermaster.js

[Install]
WantedBy=multi-user.target
