# onlyoffice-fileconverter.service
[Unit]
Description=Onlyoffice documentserver
After=network.target

[Service]
User=onlyoffice
Type=simple
Environment="NODE_ENV=production-linux" "NODE_CONFIG_DIR=/etc/webapps/onlyoffice/documentserver" "LD_LIBRARY_PATH=/usr/share/webapps/onlyoffice/documentserver/server/FileConverter/bin"
#WorkingDirectory=/usr/share/webapps/onlyoffice/documentserver/server/FileConverter
ExecStart=/usr/share/webapps/onlyoffice/documentserver/server/FileConverter/converter
WorkingDirectory=/var/lib/onlyoffice/documentserver
LimitNOFILE=1048576
LimitNPROC=64
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/onlyoffice/documentserver

[Install]
WantedBy=multi-user.target
