# oor.service
[Unit]
Description=Open Overlay Router LISP Daemon

[Service]
ExecStart=/usr/bin/oor

[Install]
WantedBy=multi-user.target
