# openhab2.service
[Unit]
Description=OpenHAB2 home automation software
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/opt/openhab2/start.sh
ExecReload=/usr/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
