# opensimulator.service
[Unit]
Description=Opensimulator
After=local-fs.target network.target

[Service]
Type=simple
ExecStart=/usr/bin/opensimulator
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always
User=%I

[Install]
WantedBy=multi-user.target

