# opennebula.service
[Unit]
Description=OpenNebula management

[Service]
User=oneadmin
Type=oneshot
ExecStart=/srv/cloud/one/bin/one start
ExecStop=/srv/cloud/one/bin/one stop
TimeoutSec=0
StandardInput=tty
RemainAfterExit=yes
