# opennebula-oneflow.service
[Unit]
Description=OpenNebula econe Server
After=network.target local-fs.target remote-fs.target opennebula.service
Requires=opennebula.service

[Service]
User=oneadmin
Group=cloud
Type=forking
PermissionsStartOnly=true
ExecStart=/usr/bin/oneflow-server start
ExecStop=/usr/bin/oneflow-server stop
Restart=on-abnormal

[Install]
WantedBy=multi-user.target
