# openbazaard-standalone.service
[Unit]
Description=Openbazaar Marketplace Server
After=network.target

[Service]
Type=forking
User=openbazaar
WorkingDirectory=/var/lib/openbazaard-standalone/
TimeoutSec=0
EnvironmentFile=/etc/conf.d/openbazaard-standalone.conf
ExecStart=/var/lib/openbazaard-standalone/openbazaard-standalone start $ARGS
ExecStop=/var/lib/openbazaard-standalone/openbazaard-standalone stop

## Security
InaccessibleDirectories=/home

[Install]
WantedBy=multi-user.target

