# openbazaard.service
[Unit]
Description=Openbazaar Marketplace Server
After=network.target

[Service]
Type=simple
User=openbazaar
DynamicUser=yes
EnvironmentFile=/etc/conf.d/openbazaard
StateDirectory=openbazaard
ExecStart=/usr/bin/openbazaard start $OB_ARGS
ExecStop=/usr/bin/openbazaard stop
ExecReload=/usr/bin/openbazaard restart

[Install]
WantedBy=multi-user.target

