# openbazaard.service
[Unit]
Description=Openbazaar Marketplace Server
After=network.target

[Service]
Type=simple
User=openbazaar
DynamicUser=yes
EnvironmentFile=/etc/conf.d/openbazaard
StateDirectory=openbazaar
ExecStart=/usr/bin/openbazaard start $OB_ARGS

[Install]
WantedBy=multi-user.target

