# opennebula-econe.service
[Unit]
Description=OpenNebula econe Server
After=syslog.target network.target remote-fs.target opennebula.service
BindsTo=opennebula.service

[Service]
User=oneadmin
Group=cloud
Type=simple

ExecStart=/usr/bin/ruby /usr/lib/one/ruby/cloud/econe/econe-server.rb
PIDFile=/var/run/one/econe-server.pid
Restart=always

[Install]
WantedBy=multi-user.target
