# openssl_version_workaround.service
[Service]
Environment="LD_LIBRARY_PATH=/opt/mssql/lib"
