# opennebula-novnc.service
[Unit]
Description=OpenNebula noVNC Server
After=syslog.target network.target remote-fs.target
Before=opennebula-sunstone.service

[Service]
User=oneadmin
Group=cloud
Type=forking

ExecStart=/usr/bin/novnc-server start
PIDFile=/var/lock/one/.novnc.lock
Restart=always

[Install]
WantedBy=multi-user.target
