# opennebula.service
[Unit]
Description=OpenNebula management

[Service]
User=oneadmin
Group=cloud
Type=forking
ExecStart=/usr/bin/one start
ExecStop=/usr/bin/one stop

[Install]
WantedBy=multi-user.target
