# opennebula-sunstone.service
[Unit]
Description=OpenNebula Web Interface
After=syslog.target network.target remote-fs.target opennebula.service opennebula-novnc.service
BindsTo=opennebula.service
Wants=opennebula-novnc.service

[Service]
User=oneadmin
Group=cloud
Type=simple
PermissionsStartOnly=true

ExecStart=/usr/bin/ruby /usr/lib/one/sunstone/sunstone-server.rb
PIDFile=/var/run/one/sunstone.pid
Restart=always

[Install]
WantedBy=multi-user.target
