# opennebula-oneflow.service
[Unit]
Description=OpenNebula Flow Service
After=syslog.target network.target remote-fs.target opennebula.service
BindsTo=opennebula.service

[Service]
User=oneadmin
Group=cloud
Type=simple


ExecStart=/usr/bin/ruby /usr/lib/one/oneflow/oneflow-server.rb
PIDFile=/var/run/one/oneflow.pid
Restart=always

[Install]
WantedBy=multi-user.target
