# openmeetings.service
# http://solie.co.vu/openmeetings-systemd-unit/
# Openmeetings Systemd Unit
# Written by Maz Solie the 18 Apr 2015 - Tagged: systemd | arch | linux | virtualbox | openmeetings

# This was originally spelled openmeeting.service (forgot the s)

[Unit]
Description=Apache OpenMeetings
After=network.target

[Service]
Type=simple
ExecStart=/opt/openmeetings/red5.sh
ExecStop=/opt/openmeetings/red5-shutdown.sh
RestartSec=5
Restart=on-failure

[Install]
WantedBy=multi-user.target
