# openbazaar.service
[Unit]
Description=Openbazaar Marketplace

[Service]
Type=forking
User=openbazaar
WorkingDirectory=/var/lib/openbazaar/
ExecStart=/usr/bin/openbazaar start
ExecStop=/usr/bin/openbazaar stop

[Install]
WantedBy=multi-user.target
