# openbazaard-standalone.service
[Unit]
Description=Openbazaar Marketplace Server
After=network.target

[Service]
Type=forking
User=openbazaar
WorkingDirectory=/var/lib/openbazaard/
TimeoutSec=0
EnvironmentFile=/etc/conf.d/openbazaard.conf
ExecStart=/var/lib/openbazaard/openbazaard start $ARGS
ExecStop=/var/lib/openbazaard/openbazaard stop
PIDFile=/tmp/openbazaard.pid

## Security
InaccessibleDirectories=/home

[Install]
WantedBy=multi-user.target

