# openbazaar-server.service
[Unit]
Description=Openbazaar Marketplace Server

[Service]
Type=forking
User=openbazaar-server
WorkingDirectory=/var/lib/openbazaar-server/
ExecStart=/usr/bin/openbazaar-server start
ExecStop=/usr/bin/openbazaar-server stop
#PIDFile=/tmp/openbazaard.pid

[Install]
WantedBy=multi-user.target
