# opennebula-onegate.service
[Unit]
Description=OpenNebula onegate Server
After=network.target local-fs.target remote-fs.target opennebula.service
Requires=opennebula.service

[Service]
User=oneadmin
Group=cloud
Type=forking
PermissionsStartOnly=true
ExecStart=/usr/bin/onegate-server start
ExecStop=/usr/bin/onegate-server stop
Restart=on-abnormal

[Install]
WantedBy=multi-user.target
