# openhab.service
[Unit]
Description=OpenHAB home automation software
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStart=/opt/openhab/start.sh server
ExecReload=/usr/bin/kill -HUP $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target
