# opentabletdriver.service
[Unit]
Description=OpenTabletDriver Daemon
After=graphical.target
PartOf=graphical.target
StartLimitIntervalSec=60
StartLimitBurst=2

[Service]
ExecStart=/usr/bin/otd
Restart=always
RestartSec=10s

[Install]
WantedBy=graphical.target
