# opentracker-ipv6.service
[Unit]
Description=Opentracker IPv6 service
After=network.target

[Service]
ExecStart=/usr/bin/opentracker-ipv6 -f /etc/opentracker-ipv6/opentracker-ipv6.conf -u otracker6
Type=simple
Restart=on-failure

[Install]
WantedBy=multi-user.target

