# opennebula-sunstone.service
[Unit]
Description=OpenNebula Web Interface
After=network.target local-fs.target remote-fs.target opennebula.service
Requires=opennebula.service

[Service]
User=oneadmin
Group=cloud
Type=forking
PermissionsStartOnly=true
ExecStart=/usr/bin/sunstone-server start
ExecStop=/usr/bin/sunstone-server stop
Restart=on-abnormal

[Install]
WantedBy=multi-user.target
