# oragono.service
[Unit]
Description=A modern IRC server written in Go
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/oragono run --conf /etc/oragono.conf
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target
