# oragono.service
[Unit]
Description=A modern IRC server written in Go

[Service]
WorkingDirectory=/opt/oragono
Type=simple
ExecStart=/opt/oragono/oragono run

[Install]
WantedBy=multi-user.target
