# orphan-manager.service
[Unit]
Description=Removes orphaned packages

[Service]
Type=simple
ExecStart=/usr/local/bin/orphans
