# orphan-manager.service
[Unit]
Description=Remove orphaned packages

[Service]
Type=simple
ExecStart=/usr/local/bin/orphan-manager
