# ory-oathkeeper@.service
[Unit]
Description=Ory Oatheeper identity and access proxy written in Go (using %i.conf)
Documentation=https://github.com/ory/oathkeeper
After=network.target

[Service]
User=ory-oathkeeeper
Group=ory-oathkeeper
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectHostname=true
ProtectClock=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
NoNewPrivileges=true
PrivateDevices=true
DeviceAllow=/dev/syslog
RestrictSUIDSGID=true
ProtectKernelModules=true
MemoryDenyWriteExecute=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
TimeoutStopSec=7s
KillMode=mixed
KillSignal=SIGTERM
Restart=on-failure
RestartPreventExitStatus=2
ExecStart=/usr/bin/ory-oathkeeper --config /etc/ory/oathkeeper/%i.yaml

[Install]
WantedBy=multi-user.target
