# osqueryd.service
[Unit]
Description=The osquery Daemon
After=network.service syslog.service

[Service]
TimeoutStartSec=0
ExecStart=/usr/bin/osqueryd --config_path /etc/osquery/osquery.conf
Restart=on-failure
KillMode=process
KillSignal=SIGTERM
SuccessExitStatus=143
PIDFile=/var/osquery/osqueryd.pidfile

[Install]
WantedBy=multi-user.target
