# ossec-hids.service
[Unit]
Description=Open Source Host-based Intrusion Detection System

[Service]
Type=forking
#User=ossec
#Group=ossec
ExecStart=/var/lib/ossec-hids/bin/ossec-control start
ExecStop=/var/lib/ossec-hids/bin/ossec-control stop

[Install]
WantedBy=basic.target
