# overture.service
[Unit]
Description=A DNS upstream switcher written in Go in order to purify DNS records
After=network.target

[Service]
ExecStart=/usr/bin/overture -c /etc/overture/config.json
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
