# owncloud-news-updater.timer
[Unit]
Description=OwnCloud News Updater Timer

[Timer]
OnBootSec=15min
OnUnitInactiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
