# p-rout-dump.timer
[Unit]
Description=Dump p-rout SQL data

[Timer]
OnCalendar=*-*-* 00:00:00
AccuracySec=1h
Persistent=true

[Install]
WantedBy=basic.target
