# paccache-leave3.service
[Unit]
Description=Clear pacman cache, leave up to 3 versions of each installed package
RequiresMountsFor=/var/cache/pacman/pkg/

[Service]
Type=oneshot
ExecStart=/usr/bin/sh -c 'paccache -vruk0 && paccache -rv'
Nice=19
StandardOutput=journal
StandardError=journal
