# paccache-leave3.timer
[Unit]
Description=Run	paccache-leave3.service weekly

[Timer]
OnCalendar=weekly
AccuracySec=24h
Persistent=yes

[Install]
WantedBy=multi-user.target
