# pacgui.timer
[Unit]
Description=Runs PACGUI daily

[Timer]
OnCalendar=daily
Persistent=true
Unit=pacgui.service

[Install]
WantedBy=default.target
