# packetbeat.service
[Unit]
Description=packetbeat
Documentation=https://www.elastic.co/guide/en/beats/packetbeat/current
After=network.target

[Service]
ExecStart=/usr/bin/packetbeat -path.config /etc/packetbeat \
                              -path.home /usr/share/packetbeat \
                              -path.data /var/lib/packetbeat \
                              -path.logs /var/log/packetbeat

[Install]
WantedBy=multi-user.target
