# pacmon.service
[Unit]
Description=pacmon
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/usr/bin/pacmon
User=pacmon
Restart=always

[Install]
WantedBy=multi-user.target
