# pacman-auto-update.timer
[Unit]
Description=Automatic Update when booted up after 5 minutes then check the system for updates every 60 minutes
[Timer]
OnBootSec=5min
OnUnitActiveSec=60min
Unit=pacman-auto-update.service
[Install]
WantedBy=multi-user.target
