# panopta-agent.cron
* * * * *    panopta-agent    /usr/bin/python /usr/bin/panopta-agent/panopta_agent.py --from-cron &> /dev/null
