# paperless-scheduler.service
[Unit]
Description=Paperless scheduler service
After=network.target
Requires=redis.service
PartOf=paperless.target

[Service]
User=paperless
Group=paperless
WorkingDirectory=/usr/share/paperless/src
ExecStart=/usr/bin/paperless-manage qcluster
Restart=on-abort

[Install]
WantedBy=paperless.target
