# paperless-consumer.service
[Unit]
Description=Paperless consumer

[Service]
User=http
Group=http
ExecStart=python /usr/share/webapps/paperless/src/manage.py document_consumer

[Install]
WantedBy=multi-user.target

