# paperless-task-queue.service
[Unit]
Description=Paperless Celery Workers
After=network.target
Requires=redis.service
PartOf=paperless.target

[Service]
User=paperless
Group=paperless
WorkingDirectory=/usr/share/paperless/src
ExecStart=/usr/bin/celery --app paperless worker --loglevel INFO
Restart=on-abort

[Install]
WantedBy=paperless.target