# parprouted.service
[Unit]
Description=Proxy-ARP bridging daemon

[Service]
Type=forking
EnvironmentFile=/etc/conf.d/parprouted
ExecStart=/usr/bin/parprouted $EXTIF $INTIF

[Install]
WantedBy=multi-user.target
