# parprouted.service
[Unit]
Description=proxy ARP IP bridging daemon

[Install]
WantedBy=multi-user.target

[Service]
Type=forking
EnvironmentFile=/etc/conf.d/parprouted
ExecStart=/usr/bin/parprouted $EXTIF $INTIF

