# paste-binoused.service
[Unit]
Description=Paste Binouse server

[Service]
Type=simple
EnvironmentFile=/etc/conf.d/paste-binoused
SyslogIdentifier=paste-binoused
ExecStart=/usr/bin/paste-binoused $PASTE_BINOUSE_ARGS
CPUSchedulingPolicy=batch

[Install]
WantedBy=multi-user.target
