# paxctld.service
[Unit]
Description=PaX flags maintenance daemon (paxctld)

[Service]
Type=forking
ExecStart=/usr/bin/paxctld -d

[Install]
WantedBy=multi-user.target
