# peerflix-server.service
[Unit]
Description=Peerflix Server

[Service]
Type=forking
User=root
WorkingDirectory=/usr/lib/node_modules/peerflix-server/
ExecStart=/usr/bin/peerflix-server

[Install]
WantedBy=multi-user.target
