# peerunityd@.service
[Unit]
Description=Peerunity Daemon
After=network.target 

[Service]
User=%i
Group=users
Type=simple
ExecStart=/usr/bin/peerunityd

[Install]
WantedBy=multi-user.target
