# peerunity-tor@.service
[Unit]
Description=Peerunity
After=network.target
After=graphical.target
Requires=tor.service

[Service]
User=%i
Group=users
Type=simple
Restart=on-failure
RestartSec=15
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
ExecStart=/usr/bin/peerunity -proxy=127.0.0.1:9050
## security
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
NoNewPrivileges=true

[Install]
WantedBy=default.target
