# peerflix-server.service
[Unit]
Description=Peerflix Server

[Service]
Type=forking
User=peerflix
ExecStart=/usr/bin/peerflix-server

[Install]
WantedBy=multi-user.target
