# peercoind@.service
[Unit]
Description=Peercoin Daemon
After=network.target 

[Service]
User=%i
Group=users
Type=simple
Restart=on-failure
RestartSec=15
ExecStart=/usr/bin/peercoind
## security
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
NoNewPrivileges=true

[Install]
WantedBy=default.target
