# pfclient-beta.service
[Unit]
Description=Plane Finder ADS-B Client
Wants=network-online.target
After=syslog.target network-online.target

[Service]
Restart=always
RestartSec=30
ExecStart=/usr/bin/pfclient-beta --daemon --log_path=/var/log --config_path=/etc/pfclient/pfclient-config.json --pid_file=/var/run/pfclient-beta.pid
ExecStop=/usr/bin/kill -INT $MAINPID
Type=forking
PIDFile=/var/run/pfclient-beta.pid

[Install]
WantedBy=multi-user.target
