# pgweb.service
[Unit]
Description=Web client (GUI) for PostgreSQL databases
After=syslog.target network.target

[Service]
Type=simple
EnvironmentFile=/etc/conf.d/pgweb
ExecStart=/usr/bin/pgweb $PGWEB_OPTIONS
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivilegies=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ReadOnlyPaths=/

[Install]
WantedBy=multi-user.target
