# phantomsocks.service
[Unit]
Description=Phantomsocks process
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/usr/bin/phantomsocks-init
Restart=no
#RestartSec=5

[Install]
WantedBy=default.target
