# phantomsocks@.service
[Unit]
Description=Phantomsocks process, on %I
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
Group=phantomsocks
WorkingDirectory=/etc/phantomsocks/%I
ExecStart=/usr/bin/phantomsocks
Restart=on-abnormal
RestartSec=5

[Install]
WantedBy=multi-user.target
Alias=
