# photoprism.service
# Do not make manual changes to this file.
# It will be overridden by package upgrades.
# Use `systemctl edit photoprism` instead!

[Unit]
After=network.target
Description=Photoprism

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=/usr/bin/photoprism --defaults-yaml=/etc/photoprism/options.yml start
Restart=always
RestartSec=10

User=photoprism
