# php-fpm.service
[Unit]
Description=The PHP 5.4 FastCGI Process Manager
After=syslog.target network.target

[Service]
Type=notify
PIDFile=/run/php54-fpm/php-fpm.pid
PrivateTmp=true
ExecStart=/usr/bin/php54-fpm --nodaemonize --pid /run/php54-fpm/php-fpm.pid
ExecReload=/bin/kill -USR2 $MAINPID

[Install]
WantedBy=multi-user.target
