# php-fpm-aegir.service
[Unit]
Description=The PHP FastCGI Process Manager
After=syslog.target network.target
Conflicts=php-fpm.service

[Service]
Type=notify
PIDFile=/run/php-fpm/php-fpm.pid
PrivateTmp=true
ExecStart=/usr/bin/php-fpm --nodaemonize --pid /run/php-fpm/php-fpm.pid
ExecReload=/bin/kill -USR2 $MAINPID

[Install]
RequiredBy=aegir.target
