# pi-hole-whitelist.timer
[Unit]
Description=Weekly update whitelist of pi-hole

[Timer]
OnCalendar=weekly
AccuracySec=3h
Persistent=true