# pi-hole-weekly.cron
#!/bin/bash
/usr/bin/gravity.sh > /dev/null

