# pi-hole-ftl.service
[Unit]
Description=Pi-hole FTL engine

[Service]
User=pihole
Group=pihole
PIDFile=/run/pihole-ftl/pihole-FTL.pid
ExecStart=/usr/bin/pihole-FTL no-daemon
Restart=on-failure
     
[Install]
WantedBy=multi-user.target
