# pi-hole-gravity.service
[Unit]
Description=Ad-serving domains collector

[Service]
Type=oneshot
ExecStart=/usr/bin/gravity.sh
Nice=19
IOSchedulingClass=best-effort
IOSchedulingPriority=7
