# pi-hole-whitelist.service
[Unit]
Description=Update whitelist of pi-hole

[Service]
Type=oneshot
ExecStart=/opt/pi-hole-whitelist/whitelist.py
