# pi-hole-ftl.service
[Unit]
Description=Pi-hole FTL engine

[Service]
Type=forking
User=pihole
Group=pihole
ExecStart=/usr/bin/pihole-FTL
Restart=on-failure
     
[Install]
WantedBy=multi-user.target
