# pi-hole-logtruncate.service
[Unit]
Description=Reset dnsmasq/pi-hole query log

[Service]
Type=oneshot
ExecStart=/usr/bin/truncate -s 0 /run/log/pihole/pihole.log
Nice=19
IOSchedulingClass=best-effort
IOSchedulingPriority=7
