# pi-hole-ftl.service
[Unit]
Description=Pi-hole FTL engine

[Service]
User=pihole
Group=pihole
ExecStart=/usr/bin/pihole-FTL
     
[Install]
WantedBy=multi-user.target
