# pi-hole-daily.cron
#!/bin/bash
truncate -s 0 /run/log/pihole.log > /dev/null

