# pi-launcher@.service
# Forked off xlogin@.service

[Unit]
Description=Qt-compositor login for user %i
After=x@vt7.service systemd-user-sessions.service
Wants=x@vt7.service
Conflicts=getty@tty7.service

[Service]
User=%i
WorkingDirectory=~
TTYPath=/dev/tty7
PAMName=login
ExecStart=-/opt/pi-compositor/pi-compositor -platform eglfs -plugin CECKeyboard

[Install]
WantedBy=graphical.target
