# pia-tools-update.timer
[Unit]
Description=Weekly Update of the pia-tools config files

[Timer]
OnCalendar=weekly
Persistent=true

