# pia@.service
[Unit]
Description=PIA OpenVPN connection to %i
After=network.target

[Service]
Type=forking
# EnvironmentFile=/etc/pia-tools.conf
ExecStart=/usr/sbin/openvpn --auth-nocache --cd /etc/openvpn/pia --config /etc/openvpn/pia/%i.ovpn --daemon openvpn@%i

[Install]
WantedBy=multi-user.target
