# pilight.service
[Unit]
Description=pilight daemon
After=network.target

[Service]
ExecStart=/usr/bin/pilight-daemon -D

[Install]
WantedBy=multi-user.target
