# pixelfed-scheduler.timer
[Unit]
Description=Pixelfed scheduler

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
