# pkg.install
post_install() {
    useradd -d /home/flow -m -r -s /usr/bin/nologin flow
    chown -R flow:flow /home/flow
}

post_upgrade() {
    systemctl restart flow-causal
}

pre_remove() {
    systemctl stop flow-causal
    systemctl disable flow-causal
    userdel flow
}
