# pkgproxy.service
[Unit]
Description=pkgproxy
After=network.target

[Service]
ExecStart=/usr/bin/pkgproxy

[Install]
WantedBy=default.target

