# plex-email.timer
[Unit]
Description=Run plex-email weekly

[Timer]
OnCalendar=Fri 18:00
Persistent=true     
 
[Install]
WantedBy=timers.target
