# plex-email.service
[Unit]
Description=plex-email Service

[Service]
GuessMainPID=no
Type=oneshot

PermissionsStartOnly=true
ExecStart=/usr/bin/python2 /opt/plex-email/plexEmail.py -c /etc/plex-email/config.conf

[Install]
WantedBy=default.target
