# plexmediaplayer@.service
[Unit]
Description = Plex Media Player
After = systemd-user-sessions.service
Conflicts = getty@tty1.service

[Service]
User = %i
PAMName = login
Type = simple
TTYPath = /dev/tty1
ExecStart = /usr/bin/xinit /usr/bin/plexmediaplayer --tv --fullscreen --scale-factor 1 -- :0 -nolisten tcp vt1
Restart = on-success
StandardInput = tty

[Install]
WantedBy = multi-user.target
