# plexrequests.service
[Unit]
Description=PlexRequests
After=network.target

[Service]
User=plexrequests
Group=plexrequests
Type=simple
WorkingDirectory=/opt/plexrequests
ExecStart=/usr/bin/meteor
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target
