# plume.service
[Unit]
Description=plume

[Service]
Type=simple
User=plume
WorkingDirectory=/usr/share/webapps/plume/
ExecStart=/usr/bin/plume
TimeoutSec=30
Restart=always

[Install]
WantedBy=multi-user.target

