# polisd@.service
[Unit]
Description=Polis Daemon
After=network.target

[Service]
User=%i
Group=users
Type=simple
ExecStart=/usr/bin/polisd

[Install]
WantedBy=multi-user.target
