# pommed.service
[Unit]
Description=Pommed Apple Hotkeys Daemon

[Service]
User=root
ExecStart=/usr/bin/pommed -f

[Install]
WantedBy=multi-user.target

