# pommed.service
[Unit]
Description=Apple laptops hotkeys events handler
Documentation=man:pommed(1)

[Service]
ExecStart=/usr/bin/pommed -f

[Install]
WantedBy=multi-user.target
