# portainer.service
[Unit]
Description=Portainer

[Service]
Type=simple
Restart=always
RestartSec=3
Environment=BIND=:9000
Environment=TUNNEL_PORT=8000
ExecStart=/usr/bin/portainer --bind $BIND --tunnel-port $TUNNEL_PORT --data /var/lib/portainer --assets /usr/share/portainer

[Install]
WantedBy=multi-user.target
