# postgresql.service
[Unit]
Description=PostgreSQL-1C database server
After=network.target

[Service]
Type=simple
TimeoutSec=120
User=postgres
Group=postgres

Environment=PGROOT=/var/lib/postgres

SyslogIdentifier=postgres
PIDFile=/var/lib/postgres/postmaster.pid
ExecStart=/usr/bin/postgres -D ${PGROOT}/data
PrivateTmp=yes
OOMScoreAdjust=-200

[Install]
WantedBy=multi-user.target

