# power.sh.service
[Service]
ExecStart=/usr/bin/sh -c '/usr/bin/udevadm settle; /usr/share/power.sh 1'
ExecStop=/usr/share/power.sh 0
RemainAfterExit=yes
