# power.sh.rules
SUBSYSTEM=="usb", ACTION=="add", ATTR{power/wakeup}=="disabled", ATTR{power/wakeup}="enabled"
SUBSYSTEM=="power_supply", ATTR{online}=="0", RUN+="/usr/bin/systemctl start power.sh.service"
SUBSYSTEM=="power_supply", ATTR{online}=="1", RUN+="/usr/bin/systemctl stop power.sh.service"
