# pppoecheck.timer
[Unit]
Description=PPPOE Check Service

[Timer]
OnActiveSec=0s
OnUnitActiveSec=60s
Persistent=true

[Install]
WantedBy=timers.target

