# prelink.timer
[Unit]
Description=Daily prelinking updates

[Timer]
OnCalendar=daily
AccuracySec=2h
Persistent=true
