# prelude-correlator.service
[Unit]
Description=Correlator of events received by Prelude

[Service]
ExecStart=/usr/sbin/prelude-correlator

[Install]
WantedBy=multi-user.target
