# prest.service
[Unit]
Description=RESTful API Server for PostgreSQL
Documentation=https://postgres.rest/
After=network-online.target postgresql.service docker.service

[Service]
Restart=on-failure
ExecStart=/usr/bin/prest
WorkingDirectory=/etc
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
