# pretix-periodic.timer
[Unit]
Description=Run pretix periodic

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target

