# prometheus-node-exporter.service
[Unit]
Description=Prometheus exporter for machine metrics
After=network.target

[Service]
EnvironmentFile=-/etc/conf.d/prometheus_node_exporter
ExecStart=/usr/bin/prometheus_node_exporter $NODE_EXPORTER_OPTS
User=nobody

[Install]
WantedBy=multi-user.target
