# prometheus-nginxlog-exporter.service
[Unit]
Description=Prometheus NginxLog Exporter
After=network.target

[Service]
ExecStart=/usr/bin/prometheus-nginxlog-exporter -config-file "/etc/prometheus/nginxlog.yml"
Restart=always
ProtectSystem=full
CapabilityBoundingSet=


[Install]
WantedBy=multi-user.target
