# prometheus-node-exporter.service
[Unit]
Description=Prometheus exporter for machine metrics

[Service]
ExecStart=/usr/bin/prometheus_node_exporter

[Install]
WantedBy=default.target
