# prometheus-node-exporter.service
[Unit]
Description=Prometheus exporter for machine metrics
After=network.target

[Service]
ExecStart=/usr/bin/prometheus_node_exporter
User=nobody

[Install]
WantedBy=multi-user.target
