# prometheus-blackbox-exporter.service
[Unit]
Description=Prometheus Exporter for machine metrics
After=network.target

[Service]
ExecStart=/usr/bin/prometheus_blackbox_exporter

[Install]
WantedBy=multi-user.target
