# prometheus-snmp-exporter.service
[Unit]
Description=Prometheus Exporter for SNMP metrics
After=network.target

[Service]
ExecStart=/usr/bin/prometheus_snmp_exporter

[Install]
WantedBy=multi-user.target