# prometheus-wireguard-exporter.service
[Unit]
Description=Prometheus WireGuard Exporter
Wants=network-online.target
After=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/conf.d/prometheus-wireguard-exporter
ExecStart=/usr/bin/prometheus_wireguard_exporter $WIREGUARD_EXPORTER_ARGS

[Install]
WantedBy=multi-user.target
