# prometheus-pushgateway.service
[Unit]
Description=Prometheus push gateway
After=network.target

[Service]
ExecStart=/usr/bin/prometheus-pushgateway

[Install]
WantedBy=multi-user.target
