# prometheus-node-exporter.service
[Unit]
Description=Prometheus Exporter for machine metrics
After=network.target

[Service]
ExecStart=/usr/bin/prometheus_node_exporter

[Install]
WantedBy=multi-user.target
