# prometheus-adguard-exporter.service
[Unit]
Description=Prometheus Exporter for AdguardHome metrics
Requires=network-online.target
After=network-online.target

[Service]
User=adguard_exporter
Group=adguard_exporter
Restart=on-failure
EnvironmentFile=-/etc/conf.d/prometheus-adguard-exporter
ExecStart=/usr/bin/prometheus-adguard-exporter $ADGUARD_EXPORTER_ARGS
ExecReload=/bin/kill -HUP $MAINPID
NoNewPrivileges=true
ProtectHome=read-only
ProtectSystem=strict

[Install]
WantedBy=multi-user.target