# prowlarr.service
# Be sure to use `systemctl edit prowlarr` to modify this service with an override.conf because
# direct changes will be overwritten by package updates.

[Unit]
Description=prowlarr Service
Wants=network-online.target
After=network-online.target

[Service]
User=prowlarr
Group=prowlarr
SyslogIdentifier=prowlarr
ExecStart=/usr/lib/prowlarr/bin/Prowlarr -nobrowser -data=/var/lib/prowlarr
Type=simple

[Install]
WantedBy=multi-user.target
