# prowlarr.service
[Unit]
Description=prowlarr Service
Wants=network-online.target
After=network-online.target

[Service]
User=prowlarr
Group=prowlarr
SyslogIdentifier=prowlarr
ExecStart=/usr/lib/prowlarr/bin/Prowlarr -nobrowser -data=/var/lib/prowlarr
Type=simple

[Install]
WantedBy=multi-user.target
