# proxy.service
# -*- coding: utf-8 -*-
# region header
# Copyright Torben Sickert (t.sickert["~at~"]gmail.com) 16.12.2012

# License
# -------

# This library written by Torben Sickert stand under a creative commons naming
# 3.0 unported license. see http://creativecommons.org/licenses/by/3.0/deed.de
# endregion
[Unit]
Description=HTTP proxy server.
After=docker.service
BindsTo=docker.service
Conflicts=shutdown.target reboot.target halt.target

[Service]
User=service
Restart=always
RestartSec=10
WorkingDirectory=/srv/http/proxy
ExecStartPre=-/usr/bin/env docker-compose --file base.yml kill $APP
ExecStartPre=-/usr/bin/env docker-compose --file base.yml rm --force $APP
ExecStart=/usr/bin/env docker-compose --file base.yml up --no-build --force-recreate --no-deps $APP
ExecStop=/usr/bin/env docker-compose --file base.yml stop $APP
NotifyAccess=all

[Install]
WantedBy=multi-user.target
# region vim modline
# vim: set tabstop=4 shiftwidth=4 expandtab:
# vim: foldmethod=marker foldmarker=region,endregion:
# endregion
