# prysm-beacon-chain.service
[Unit]
Description=Prysm Beacon chain daemon
After=network-online.target

[Service]
DynamicUser=true
ExecStart=/usr/bin/prysm-beacon-chain --accept-terms-of-use
Restart=on-failure

[Install]
WantedBy=default.target