# psiphon-server.service
[Unit]
Description=Psiphon Proxy Server Service
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=%E/psiphon/
ExecStart=/usr/bin/psiphond run
Restart=on-failure

[Install]
WantedBy=multi-user.target