# pterodactyl-queue.service
[Unit]
Description=Pterodactyl queue service
After=redis.service

[Service]
User=http
Restart=always
ExecStart=/usr/bin/php /usr/share/webapps/pterodactyl/artisan queue:work --queue=high,standard,low --sleep=3 --tries=3
StartLimitInterval=180
StartLimitBurst=30
RestartSec=5s

[Install]
WantedBy=multi-user.target
