# pulledpork.cron
#!/bin/sh

/usr/bin/pulledpork_update.sh /etc/snort/rules/snort.rules > /dev/null 2>&1
