# pump.io.service
[Unit]
Description=Pump.io Social Server
After=network.target

[Service]
Type=forking
User=pump
Group=pump
ExecStart=/usr/bin/pump.io start
ExecStop=/usr/bin/pump.io stop
Restart=always
SyslogIdentifier=pump.io

[Install]
WantedBy=multi-user.target
