# pushgateway.service
[Unit]
Description=Prometheus Pushgateway
After=network.target

[Service]
ExecStart=/usr/bin/pushgateway

[Install]
WantedBy=multi-user.target
