# pyca.target
[Unit]
Description=Python Capture Agent
Documentation=https://github.com/opencast/pyCA
Requires=multi-user.target
Wants=network.target
