# pyca-capture.service
[Unit]
Description=Python Capture Agent capture service
Documentation=https://github.com/opencast/pyCA
Wants=network.target
PartOf=pyca.service
After=pyca.service

[Service]
Type=notify
NotifyAccess=all
WatchdogSec=300
User=pyca
ExecStart=/usr/bin/pyca capture
Restart=always
RestartSec=10
TimeoutSec=300

[Install]
WantedBy=pyca.service
