# pyca.service
[Unit]
Description=Python Capture Agent
Documentation=https://github.com/opencast/pyCA
Wants=network.target
PartOf=pyca.target

[Service]
Type=simple
User=pyca
ExecStart=/usr/bin/pyca
Restart=always
RestartSec=10
TimeoutSec=300

[Install]
WantedBy=pyca.target
