# pyload.service
[Unit]
Description=pyLoad Core
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/pyload --daemon --configdir=/var/lib/pyload
Restart=on-abort
KillSignal=SIGQUIT
PIDFile=/var/lib/pyload/pyload.pid
User=pyload
Group=pyload

[Install]
WantedBy=multi-user.target
