# pyload.service
[Unit]
Description=pyLoad Daemon
After=network.target

[Service]
ExecStart=/usr/bin/pyload --daemon
User=pyload
Group=pyload

[Install]
WantedBy=multi-user.target
