# quint@.service
# Forked off xlogin@.service

[Unit]
Description=quint session for user %i
After=systemd-user-sessions.service
Conflicts=getty@tty7.service

[Service]
User=%i
WorkingDirectory=~
TTYPath=/dev/tty7
PAMName=login
ExecStart=/opt/quint/quint -platform eglfs

[Install]
WantedBy=graphical.target
Alias=display-manager.service
