# radarr.service
[Unit]
Description=Radarr Service
Wants=network-online.target
After=network-online.target

[Service]
User=radarr
Group=radarr
ExecStart=/usr/lib/radarr/Radarr -nobrowser -data=/var/lib/radarr
Type=simple
TimeoutStopSec=20

# Lets built in updater work well.
Restart=on-failure
KillMode=process

[Install]
WantedBy=multi-user.target
