# radsecproxy.service
[Unit]
Description=generic RADIUS proxy
After=network.target
Conflicts=radsecproxy-keepalive.service

[Service]
Type=forking
Restart=always
ExecStartPre=/usr/bin/radsecproxy -p
ExecStart=/usr/bin/radsecproxy
ProtectSystem=full
ProtectHome=on
PrivateDevices=on
NoNewPrivileges=on


[Install]
WantedBy=multi-user.target
