# raid-check.service
[Unit]
Description=raid-check

[Service]
Type=oneshot
ExecStart=/usr/bin/raid-check
