# raumserverDaemon.service
[Unit]
Description=raumsrv-Daemon-server
After=syslog.target network.target 

[Service]
Type=forking
WorkingDirectory=/opt/raumserverDaemon/
ExecStart=/opt/raumserverDaemon/raumsrvDaemon

[Install]
WantedBy=multi-user.target
