# ravenna-alsa.service
[Unit]
Description=daemon for ALSA RAVENNA/AES67 Driver
Requires=sound.target avahi-daemon.service
After=sound.target avahi-daemon.service

[Service]
Type=simple
ExecStart=/bin/bash /opt/ravenna-alsa/ravenna_start.sh
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
