# rbldnsd-sync.timer
[Unit]
Description=RBLDNSD Synchronisation Timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target

