# rbldnsd-sync.service
[Unit]
Description=RBLDNSD Synchronisation

[Service]
Group=rbldnsd
Type=simple
Nice=19
ExecStart=/usr/lib/rbldnsd/rbldnsd-sync.sh
StandardOutput=file:/var/rbldnsd/sync.log
StandardError=inherit

[Install]
WantedBy=default.target

