# rcsc-localsign.service
[Unit]
Description=RCSC LocalSign token server

[Service]
Type=simple
User=daemon
ExecStart=/usr/bin/java -jar /opt/RCSC/LocalSign/local_webserver.jar
KillSignal=SIGKILL

[Install]
WantedBy=multi-user.target
