# rdiffweb.service
[Unit]
Description=Rdiffweb Server
Documentation=http://www.patrikdufresne.com/en/rdiffweb/
After=network.target

[Service]
ExecStart=/usr/bin/rdiffweb

[Install]
WantedBy=multi-user.target
