# rdircd.service
[Service]
ExecStart=/usr/bin/rdircd -c /var/lib/rdircd/config.ini -c /var/lib/rdircd/state.ini
Restart=always
RestartSec=5s
User=rdircd

[Install]
WantedBy=multi-user.target
