# rdm.socket
[Unit]
Description=RTags Server

[Socket]
ListenStream=%h/.rdm

[Install]
WantedBy=default.target
