# rdm.service
[Unit]
Description=RTags daemon

[Service]
Type=simple
ExecStart=/usr/bin/rdm

[Install]
WantedBy=default.target
