# readonly-cache.service
[Unit]
Description=readonly cache service

[Service]
Type=oneshot
ExecStart=/usr/bin/readonly-cache
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
