# red-discordbot.service
[Unit]
Description=Red Discord bot daemon

[Service]
WorkingDirectory=/opt/red-discordbot/
User=red
Group=red
ExecStart=/usr/bin/python launcher.py --start

[Install]
WantedBy=multi-user.target
