# redditfs.service
[Unit]
Description=redditfs
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/redditfs
WorkingDirectory=/usr/bin/
Restart=always

[Install]
WantedBy=multi-user.target
