# remctl.socket
[Unit]
Description=remctl command execution service
Documentation=man:remctld(8)
Conflicts=remctl.service

[Socket]
ListenStream=4373
Accept=yes

[Install]
WantedBy=sockets.target
