# remctl.service
[Unit]
Description=remctl command execution service
Documentation=man:remctld(8)

[Service]
ExecStart=/usr/bin/remctld -m -F
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
