# renderd-standalone.service
[Unit]
Description=renderd
Documentation=man:renderd(1)

[Service]
Type=forking
ExecStart=/usr/bin/renderd

[Install]
WantedBy=multi-user.target
