# reposilite.service
[Unit]
Description=Repository manager service

[Service]
Type=simple
User=reposilite
ExecStart=/usr/bin/java -jar /usr/share/java/reposilite/reposilite.jar -wd /var/lib/reposilite
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
