# reresolve-dns-systemd.service
[Unit]
Description=Refresh Wireguard endpoints
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/reresolve-dns-systemd.sh

[Install]
WantedBy=network-online.target

