# resolvconf.service
[Service]
Type=oneshot
ExecStart=/usr/bin/curl -4 -s --retry 1024 -o /etc/resolv.conf "https://api.opennicproject.org/geoip/?resolv&pct=99&adm=99"
