# retroarch-standalone.service
[Unit]
Description=Systemd service and user to run Retroarch in stand-alone mode

[Service]
User=retroarch
Group=retroarch
PAMName=login
Type=simple
ExecStart=/usr/bin/retroarch --appendconfig=/etc/retroarch-standalone.cfg
TimeoutStopSec=30
Restart=always

[Install]
WantedBy=multi-user.target
