# riak.service
[Unit]
Description=Distributed key/value store from Basho Technologies
After=network.target

[Service]
User=riak
Type=forking
ExecStart=/opt/riak/bin/riak start
ExecStop=/opt/riak/bin/riak stop

[Install]
WantedBy=multi-user.target
