# ricochet-irc.service
[Unit]
Description=Ricochet bridge for IRC

[Service]
ExecStart=/usr/bin/ricochet-irc --config %h/.config/ricochet-irc/
Restart=always

[Install]
WantedBy=default.target
