# rippled.service
[Unit]
Description=Ripple Peer-to-Peer Network Daemon
 
[Service]
Type=simple
ExecStart=/usr/bin/rippled --conf /etc/ripple/rippled.cfg
LimitNOFILE=8192
 
[Install]
WantedBy=multi-user.target
