# rmilter.service
[Unit]
Description=another sendmail filter

[Service]
Type=forking

PermissionsStartOnly=true

PIDFile=/run/rmilter/rmilter.pid
User=_rmilter
ExecStartPre=/usr/bin/mkdir -p /run/rmilter
ExecStartPre=/usr/bin/chown _rmilter:_rmilter /run/rmilter
ExecStart=/usr/bin/rmilter -c /etc/rmilter/rmilter.conf
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target

