# roonserver.service
[Unit]
    Description=RoonServer
    After=network.target

    [Service]
    Type=simple
    User=root
    Environment=ROON_DATAROOT=/var/roon
    Environment=ROON_ID_DIR=/var/roon
    ExecStart=/opt/RoonServer/start.sh
    Restart=on-abort

    [Install]
    WantedBy=multi-user.target