# rplayd.socket
[Unit]
Description=Rplay Socket
Conflicts=rplayd.service

[Socket]
ListenStream=5556
Accept=true

[Install]
WantedBy=sockets.target
