# rqlite@.service
[Unit]
Description=
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rqlite
Group=rqlite
ExecStart=/usr/bin/rqlited -node-id $(cat /etc/machine-id) -http-addr 0.0.0.0:4001 -raft-addr 0.0.0.0:4002 -disco-id %i /run/rqlite
Restart=always
WorkingDirectory=/run/rqlite

[Install]
WantedBy=multi-user.target
