# rsnapshot-daily@.timer
[Unit]
Description=rsnapshot daily timer

[Timer]
OnCalendar=*-*-* 5:30
Persistent=true
Unit=rsnapshot-daily@%i.service

[Install]
WantedBy=timers.target
