# rspamd.service
[Unit]
Description=rapid spam filtering system
After=nss-lookup.target

[Service]
Type=forking
PIDFile=/run/rspamd/rspamd.pid
ExecStartPre=/usr/bin/mkdir -p /run/rspamd
ExecStartPre=/usr/bin/chown _rspamd:_rspamd /run/rspamd
ExecStart=/usr/bin/rspamd -c /etc/rspamd/rspamd.conf
User=_rspamd
PermissionsStartOnly=true

[Install]
WantedBy=multi-user.target

