# rtapp.service
[Unit]
Description=Realtime Application thread system tuning

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/rtapp
