# rtc-ds1307.service
[Unit]
Description=Raspberry Pi i2C RTC DS1307 clock boot load script
After=systemd-modules-load.service


[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/systemd/scripts/rtc-ds1307-load


[Install]
WantedBy=multi-user.target
