# RTL.service
[Unit]
Description=Ride the Lightning daemon
After=network.target
Wants=lnd.service
After=lnd.service

[Service]
Type=simple
Environment="RTL_CONFIG_PATH=/etc/RTL"
Environment="RTL_LOG_PATH=/var/log/RTL"
LogsDirectory=RTL
ConfigurationDirectory=RTL
ExecStartPre=+chown bitcoin:bitcoin /etc/RTL/RTL-Config.json
ExecStart=/usr/bin/node /opt/RTL/rtl
User=bitcoin
Group=bitcoin
Restart=on-failure
TimeoutSec=120
RestartSec=30

[Install]
WantedBy=multi-user.target
