# rtorrent.service
[Unit]
Description=rTorrent System Daemon
After=network.target

[Service]
Type=simple
User=rtorrent
Group=rtorrent
ExecStart=/usr/bin/rtorrent -o system.daemon.set=true,import=/var/lib/rtorrent/rtorrent.rc
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
