# rtsp-simple-server.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/bin/rtsp-simple-server /etc/rtsp-simple-server/rtsp-simple-server.yml

[Install]
WantedBy=multi-user.target
