# rtsp-simple-server.service
[Unit]
Description=Ready-to-use RTSP / RTMP / LL-HLS server and proxy
After=network.target

[Service]
Type=simple
User=rtsp-simple-server
ExecStart=/usr/bin/rtsp-simple-server /etc/rtsp-simple-server/rtsp-simple-server.yml

[Install]
WantedBy=multi-user.target
