# rtsp-simple-server.service
[Unit]
After=network.target

[Service]
User=root
Type=simple
ExecStart=/usr/bin/rtsp-simple-server /etc/rtsp-simple-server.yml
StandardOutput=file:/run/rtsp-simple-server.log
StandardError=file:/run/rtsp-simple-server.log
ExecStop=/bin/kill -s STOP $MAINPID

[Install]
WantedBy=multi-user.target
