# runrestic.service
[Unit]
Description=runrestic backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/runrestic
