# runsunloginclient.service
[Unit]
Description=check running sunloginclinet loop
Requires=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStart=/opt/sunlogin/bin/oray_rundaemon -m server
KillMode=process
ExecStop=

[Install]
WantedBy=multi-user.target
