# rusty-keys.target
[Unit]
Description=rusty-keys keymapper service
After=systemd-udevd.service

[Install]
WantedBy=default.target

