# ryzen-ppd.service
[Unit]
Description=Ryzen Power Profile Daemon

[Service]
Type=simple
ExecStart=/usr/bin/ryzen-ppd

[Install]
WantedBy=multi-user.target
