# sabnzbd-docker-update.timer
[Unit]
Description=Run sabnzbd-docker update daily and on boot
After=docker.service

[Timer]
OnBootSec=15min
OnUnitActiveSec=1w

[Install]
WantedBy=timers.target
