# sabnzbd.service
[Unit]
Description=SABnzbd binary newsreader
After=network.target

[Service]
ExecStart=/opt/sabnzbd/SABnzbd.py -l0 -f /opt/sabnzbd/sabnzbd.ini
User=sabnzbd
Group=sabnzbd

[Install]
WantedBy=multi-user.target
