# sakura-frpc@.service
[Unit]
Description=Sakura Frp Client Service
# After=network.target
Wants=network-online.target
After=network-online.target

[Service]
Type=idle
DynamicUser=yes
Restart=on-failure
RestartSec=5s
ExecStart=/usr/bin/sakura-frpc -c /etc/sakura-frp/%i.ini
ExecReload=/usr/bin/sakura-frpc reload -c /etc/sakura-frp/%i.ini

[Install]
WantedBy=multi-user.target
