# sandpolis-server.service
[Unit]
Description=The Sandpolis Server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=sandpolis
ExecStart=/usr/bin/sandpolis-server

[Install]
WantedBy=multi-user.target