# sanoid.service
[Unit]
Description=Sanoid

[Service]
Type=oneshot
ExecStart=/usr/bin/sanoid --cron --verbose
