# sanoid.timer
[Unit]
Description=Sanoid timer

[Timer]
OnBootSec=0min
OnCalendar=minutely

[Install]
WantedBy=timers.target
