# sanoid.timer
[Unit]
Description=Sanoid timer

[Timer]
OnBootSec=0min
# Every minute, watch out for "frequently" snapshots.
#OnCalendar=minutely
# Every 15 minutes, more sane w/ "frequently" snapshots.
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
