# saturated.service
[Unit]
Description=Dead simple makepkg builder server

[Service]
User=makepkg
LimitNOFILE=65535
ExecStart=/usr/bin/saturated :8080
Restart=always

[Install]
WantedBy=multi-user.target
