# saved-kernel-modules.service
[Unit]
Description=Remove saved kernel modules on boot
ConditionPathExists=/usr/lib/modules/saved-kernel-modules

[Service]
Type=oneshot
ExecStart=/usr/lib/modules/saved-kernel-modules/cleanup
IOSchedulingClass=idle
