# scion-dispatcher.service
[Unit]
Description=SCION Dispatcher
Documentation=https://www.scionlab.org
After=network-online.target
Wants=network-online.target
PartOf=scionlab.target

[Service]
Type=simple
User=scion
Group=scion
ExecStartPre=/bin/rm -rf /dev/shm/dispatcher/
ExecStart=/usr/bin/scion-dispatcher --config /etc/scion/dispatcher.toml
RemainAfterExit=False
LimitNOFILE=4096
KillMode=control-group
Restart=on-failure
RestartSec=10

[Install]
WantedBy=scionlab.target
