# scissy.service
[Unit]
Description=Scissy git repository server

[Service]
Type=simple
ExecStart=/usr/bin/scissyd
CPUSchedulingPolicy=batch

[Install]
WantedBy=multi-user.target
