# sd-ssh-agent.socket
[Unit]
ConditionEnvironment=SSH_AUTH_SOCK

[Socket]
ListenStream=%t/sd-ssh-agent.sock
Backlog=128

[Install]
WantedBy=sockets.target
