# seafile-server@.service
[Unit]
Description=Next-generation open source cloud storage.
After=syslog.target network.target

[Service]
Type=forking
User=seafile
WorkingDirectory=/srv/seafile/%i/seafile-server

ExecStartPre=/srv/seafile/%i/seafile-server/seafile.sh start
ExecStart=/srv/seafile/%i/seafile-server/seahub.sh start
ExecStop=/srv/seafile/%i/seafile-server/seahub.sh stop
ExecStop=/srv/seafile/%i/seafile-server/seafile.sh stop

[Install]
WantedBy=multi-user.target
