# searx.service
[Unit]
Description=Searx server
After=network.target

[Service]
Type=simple
User=searx
Environment=SEARX_SETTINGS_PATH=/etc/searx/settings.yml
ExecStart=/usr/bin/searx-run
Restart=on-abort

[Install]
WantedBy=multi-user.target
