# sentry.service
[Unit]
Description=Sentry logging and aggregation server.
Requires=sentry-web.service sentry-celery.service sentry-cron.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
