# sentry-celery.service
[Unit]
Description=Sentry background worker.
After=network.target
PartOf=sentry.service

[Service]
User=sentry
Restart=always
Environment=PYTHONPATH=/opt/sentry/usr/lib/python3.10/site-packages/
ExecStart=/opt/sentry/bin/sentry --config=/etc/sentry run worker --loglevel=ERROR

[Install]
WantedBy=multi-user.target
