# sentry.service
[Unit]
Description=Sentry
After=network.target
Requires=sentry-celery.service

[Service]
User=sentry
Restart=always
Environment=PATH=/opt/sentry/bin
ExecStart=/opt/sentry/bin/sentry --config=/etc/sentry/sentry.conf.py start

[Install]
WantedBy=multi-user.target
