# sentry.target
[Unit]
Description=Sentry logging and aggregation server.
