# serviio.service
[Unit]
Description=Serviio
After=network.target

[Service]
ExecStart=/usr/bin/serviio

[Install]
WantedBy=multi-user.target
