# sfptpd.service
[Unit]
Description=sfptpd Service
After=network.target
Requires=openonload.service
After=openonload.service

[Service]
Type=forking
PrivateTmp=true
EnvironmentFile=-/etc/default/sfptpd

ExecStart=/usr/bin/sfptpd -f /etc/sfptpd.cfg
ExecStop=/usr/bin/killall -SIGTERM /usr/bin/sfptpd

[Install]
WantedBy=multi-user.target
