# sgeexecd@.service
[Unit]
Description=SGE qmaster service
After=network.target

[Service]
Type=simple
ExecStart=/opt/sge/%i/common/sgeexecd

[Install]
WantedBy=multi-user.target
