# shadowsocksr-native@.service
[Unit]
Description=ShadowsocksR-native client service for %I
After=network.target

[Service]
Type=simple
WorkingDirectory=/etc/ssr-native
ExecStart=/usr/bin/ssr-client -c %i.json
DynamicUser=yes
ProtectHome=true

[Install]
WantedBy=multi-user.target
