# shadowsocksr-libev-server@.service
[Unit]
Description=ShadowsocksR-Libev Server Service
After=network.target

[Service]
Type=simple
User=nobody
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
ExecStart=/opt/shadowsocksr-libev/bin/ssr-server -c /etc/shadowsocksr/%i.json

[Install]
WantedBy=multi-user.target

