# shellinabox@.service
[Unit]
Description=Shell In a Box Daemon

[Service]
Type=oneshot
RemainAfterExit=yes
User=%i
ExecStart=/usr/bin/shellinaboxd -b -s '/:SSH:localhost' --localhost-only -c '/tmp'

[Install]
WantedBy=multi-user.target
