# shellinabox@.service
[Unit]
Description=Shell In a Box service

[Service]
Type=oneshot
RemainAfterExit=yes
User=%i
ExecStart=shellinaboxd -b -s '/:SSH:localhost' --localhost-only -c /tmp

[Install]
WantedBy=multi-user.target
