# shiny-server.service
[Unit]
Description=Shiny Server
After=syslog.target network.target

[Service]
User=shiny
ExecStart=/usr/bin/shiny-server

[Install]
WantedBy=multi-user.target
