# shoved.service
[Unit]
Description=shoved Push Notifications Service

[Service]
ExecStart=/usr/bin/shoved
Restart=always

[Install]
WantedBy=default.target
