# sickchill.service
[Unit]
Description=SickChill Daemon
After=network.target

[Service]
Restart=on-success
User=sickchill
Group=sickchill

Type=simple
#ExecStart=/var/lib/sickchill/app/bin/SickChill --quiet --nolaunch --datadir=/var/lib/sickchill/data
ExecStart=/opt/sickchill/app/bin/SickChill --quiet --nolaunch --datadir=/opt/sickchill/data

[Install]
WantedBy=multi-user.target
