# sickrage.service
[Unit]
Description=SickRage Daemon
After=network.target

[Service]
Restart=on-success
User=sickrage
Group=sickrage

Type=forking
GuessMainPID=no
#ExecStart=/usr/bin/env python2 /opt/sickrage/app/SickBeard.py --quiet --daemon --nolaunch --datadir=/opt/sickrage/data
ExecStart=/usr/bin/python2 /opt/sickrage/app/SickBeard.py --quiet --daemon --nolaunch --datadir=/opt/sickrage/data

[Install]
WantedBy=multi-user.target
