# sickrage-user.service
[Unit]
Description=SickRage Daemon
## If you want to ensure a torrent client (or usenet client) is running)
#Requires=deluged.service

[Service]
ExecStart=/usr/bin/env python2 /opt/sickrage/SickBeard.py --quiet --config %h/.sickrage/config.ini --datadir %h/.sickrage

[Install]
WantedBy=default.target
