# sickrage-system.service
[Unit]
Description=SickRage Daemon
After=network.target

[Service]
User=sickrage
Type=forking
PIDFile=/run/sickrage/sickrage.pid
ExecStart=/usr/bin/env python2 /opt/sickrage/SickBeard.py -q --daemon --nolaunch --pidfile=/run/sickrage/sickrage.pid --config /opt/sickrage/config.ini --datadir /opt/sickrage

[Install]
WantedBy=multi-user.target
