# sickchill.service
[Unit]
Description=SickChill Daemon
After=network.target

[Service]
Restart=on-success
User=sickchill
Group=sickchill

Type=simple
#ExecStart=/usr/bin/env python2 ...
ExecStart=/usr/bin/python2 /opt/sickchill/app/SickBeard.py --quiet --nolaunch --datadir=/opt/sickchill/data

[Install]
WantedBy=multi-user.target
