# sickrage.service
[Unit]
Description=SickRage Daemon
After=network.target

[Service]
Restart=on-success
User=sickrage
Group=sickrage

#Type=forking
#PIDFile=/run/sickrage/sickrage.pid
#ExecStart=/usr/bin/env python2 /opt/sickrage/app/SickBeard.py --quiet --daemon --datadir /opt/sickrage/data --pidfile=/run/sickrage/sickrage.pid
ExecStart=/usr/bin/env python2 /opt/sickrage/app/SickBeard.py --quiet --nolaunch --datadir /opt/sickrage/data

[Install]
WantedBy=multi-user.target
