# signalk-server.socket
[Socket]
ListenStream=3443
ListenStream=3000

[Install]
WantedBy=sockets.target
