# signal-web-gateway.service
[Unit]
Description=signal-web-gateway daemon
After=network.target

[Service]
PIDFile=/var/lib/signal-web-gateway/pid
User=signal
Group=signal
RuntimeDirectory=signal-web-gateway
WorkingDirectory=/var/lib/signal-web-gateway/
ExecStart=/usr/bin/signal-web-gateway -gateway
LimitNOFILE=1048576
LimitNPROC=64
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/signal-web-gateway

[Install]
WantedBy=multi-user.target
