# simple-rt@.service
[Unit]
Description=SimpleRT Daemon on %i

[Service]
Type=idle
ExecStart=/usr/lib/simple-rt/simple-rt -i %i
WorkingDirectory=/usr/lib/simple-rt
Restart=always

[Install]
WantedBy=multi-user.target
