# simple-http-server@.service
[Unit]
Description=Simple http server in Rust
After=network.target network-online.target nss-lookup.target

[Service]
WorkingDirectory=%f
ExecStart=/usr/bin/simple-http-server -p 80

[Install]
WantedBy=multi-user.target
