# sinfod.service
[Unit]
Description=Monitoring tool for networked computers - demon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/sinfod -Fq

[Install]
WantedBy=multi-user.target
