# sinfod.service
[Unit]
Description=Monitoring tool for networked computers - demon
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/sinfod -q

[Install]
WantedBy=multi-user.target
