# sinopia.service
[Unit]
Description=Private npm cache registry

[Service]
ExecStart=/usr/bin/sinopia

[Install]
WantedBy=default.target

