# sks-recon.service
[Unit]
Description=Synchronizing key server recon instance
After=network.target remote-fs.target nss-lookup.target

[Service]
Type=simple
User=sks
Group=sks
PIDFile=/var/run/sks/sksrecon.pid
ExecStart=/usr/bin/sks recon
RuntimeDirectory=sks
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
