# sks-recon-local.service
[Unit]
Description=Synchronizing key server recon instance (Local sync/dumps)
After=network.target remote-fs.target nss-lookup.target sks-db-local.service
Requires=sks-db-local.service

[Service]
Type=simple
User=sks
Group=sks
PIDFile=/var/run/sks-local/sksrecon.pid
ExecStart=/usr/bin/sks-local recon

[Install]
WantedBy=multi-user.target
