# skywire-autoconfig.service
[Unit]
Description=skywire-autoconfig
After=systemd-networkd.service network.target armbian-firstrun.service armbian-resize-filesystem.service systemd-user-sessions.service network-online.target systend-networkd-wait-online.service NetworkManager-wait-online.service
Before=skywire.service
Wants=network-online.target

[Service]
Type=forking
Environment="SYSTEMDCHILD=1"
ExecStart=/usr/bin/skywire-autoconfig

[Install]
WantedBy=multi-user.target
