# skywire.service
[Unit]
Description=Skywire
Wants=network-online.target

[Service]
Type=simple
Environment="SYSTEMDCHILD=1"
ExecStart=/bin/skywire -p --loglvl debug'
Restart=always
RestartSec=20
TimeoutSec=30

[Install]
WantedBy=multi-user.target
