# skywire.service
[Unit]
Description=Skywire
Before=skywire-visor.service
Conflicts=skywire-visor.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
Environment="SYSTEMDCHILD=1"
ExecStart=/usr/bin/skywire -p
Restart=always
RestartSec=20
TimeoutSec=30

[Install]
WantedBy=multi-user.target
