# slskd.service
# Be sure to use `systemctl edit slskd` to modify this service with an override.conf because
# direct changes will be overwritten by package updates.
#
# A user per service w/ shared group setup would have an override like:
# [Service]
# Group=sharedgroupname
# Umask=002

[Unit]
Description=Slskd Service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=slskd
Group=slskd
SyslogIdentifier=slskd
ExecStart=/usr/lib/slskd/slskd
Restart=on-failure

[Install]
WantedBy=multi-user.target
