# slurmd.service
[Unit]
Description=Slurm node daemon
After=network.target systemd-tmpfiles-clean.service
Requires=munge.service
ConditionPathExists=/etc/slurm-llnl/slurm.conf

[Service]
Type=forking
EnvironmentFile=-/etc/default/slurm-llnl
ExecStart=/usr/bin/slurmd $SLURMD_OPTIONS
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/run/slurm-llnl/slurmd.pid
KillMode=process
LimitNOFILE=51200
LimitMEMLOCK=infinity
LimitSTACK=infinity

[Install]
WantedBy=multi-user.target
