# slurmdbd.service
[Unit]
Description=Slurm DBD accounting daemon
After=network.target mysqld.service
ConditionPathExists=/etc/slurm-llnl/slurmdbd.conf

[Service]
Type=forking
EnvironmentFile=-/etc/default/slurm-llnl
ExecStart=/usr/bin/slurmdbd $SLURMDBD_OPTIONS
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/run/slurmdbd.pid

[Install]
WantedBy=multi-user.target
