# slurmd.service
[Unit]
Description=Slurm node daemon
After=network.target systemd-tmpfiles-clean.service
Requires=munge.service
ConditionPathExists=/etc/slurm-llnl/slurm.conf

[Service]
Type=simple
EnvironmentFile=/etc/default/slurm-llnl
ExecStartPre=/usr/bin/mkdir -p /var/lib/slurm-llnl
ExecStart=/usr/bin/slurmd -D $SLURMD_OPTIONS
PIDFile=/var/run/slurm/slurmd.pid

[Install]
WantedBy=multi-user.target
