# smtpd.service
[Unit]
Description=A FREE implementation of the server-side SMTP protocol
After=network.target
Wants=network.target

[Service]
ExecStart=/usr/sbin/smtpd -d
ExecStop=/usr/bin/smtpctl stop

[Install]
WantedBy=multi-user.target
