# snort.service
[Unit]
Description=Snort IDS daemon in NFQUEUE mode

[Service]
ExecStart=/usr/bin/snort -Q -c /etc/snort/snort.lua -l /var/log/snort --tweaks local
ExecReload=kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
