# softethervpn-bridge.service
[Unit]
Description=SoftEther VPN Bridge daemon
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/vpnbridge start
ExecReload=/usr/bin/vpnbridge stop

[Install]
WantedBy=multi-user.target
