# softethervpn-client.service
[Unit]
Description=SoftEther VPN Client
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/vpnclient start
ExecStop=/usr/bin/vpnclient stop
Restart=on-failure

# Hardening
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=full
ReadOnlyDirectories=/
ReadWriteDirectories=-/usr/lib/softethervpn/vpnclient

[Install]
WantedBy=multi-user.target
