# softethervpn-bridge.service
[Unit]
Description=SoftEther VPN Bridge
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/vpnbridge start
ExecStop=/usr/bin/vpnbridge stop
Restart=on-failure

# Hardening
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=full
ReadOnlyDirectories=/
ReadWriteDirectories=-/usr/lib/softethervpn/vpnbridge

[Install]
WantedBy=multi-user.target
