# soju.service
[Unit]
Description=soju IRC bouncer service
After=network.target

[Service]
Type=simple
User=soju
Restart=always
ExecStart=/usr/bin/soju --config /etc/soju/config
ExecReload=/bin/kill -HUP $MAINPID
RuntimeDirectory=soju

[Install]
WantedBy=multi-user.target
