# solarmeter.service
[Unit]
Description=Solarmeter
After=syslog.target
#Wants=mosquitto.service

[Service]
ExecStart=/usr/bin/solarmeter -c /etc/solarmeter.conf
User=solarmeter
Group=solarmeter
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
