# someblocks.service
[Unit]
Description=Someblocks
Wants=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/someblocks

[Install]
WantedBy=default.target
