# sonarr.service
[Unit]
Description=Sonarr Service
After=network.target

[Service]
User=sonarr
Group=sonarr
ExecStart=/usr/bin/sonarr
Type=simple
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
