# sonarr.service
[Unit]
Description=Sonarr Service
After=network.target

[Service]
User=sonarr
Group=sonarr
ExecStart=/usr/bin/mono /usr/lib/sonarr/NzbDrone.exe -nobrowser -data=/var/lib/sonarr
Restart=always
Type=simple

[Install]
WantedBy=multi-user.target
