# sonarrannounced.service
[Unit]
Description=sonarrAnnounced Service

[Service]
WorkingDirectory=/opt/sonarrAnnounced
ExecStart=/usr/bin/python /usr/lib/sonarrAnnounced/bot.py
Restart=always
GuessMainPID=no
Type=simple
User=sonarrannounced
Group=sonarrannounced

[Install]
WantedBy=multi-user.target
