# sonarrannounced.service
[Unit]
Description=Bazarr Service
After=network.target

[Service]
User=sonarrannounced
Group=sonarrannounced
Type=simple
ExecStart=/usr/bin/python2 /usr/lib/sonarrannounced/sonarrannounced.py --no-update --config /var/lib/sonarrannounced
Restart=always
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
