# sonarr.service
[Unit]
Description=Sonarr Service
After=network.target

[Service]
User=sonarr
Group=sonarr
ExecStart=/usr/bin/sonarr
Type=simple
PIDFile=/var/lib/sonarr/sonarr/.config/NzbDrone/nzbdrone.pid
Restart=always

[Install]
WantedBy=multi-user.target
