# spacenavd.service
[Unit]
Description=Spacenav driver daemon

[Service]
Type=forking
PIDFile=/var/run/spnavd.pid
ExecStart=/usr/bin/spacenavd

[Install]
WantedBy=multi-user.target
