# sphinx.service
[Unit]
Description=Sphinx - SQL Full Text Search Engine
After=local-fs.target network.target mysqld.service

[Service]
Type=forking
Restart=always
PIDFile=/var/lib/sphinx/log/searchd.pid
ExecStart=/usr/bin/searchd --config /etc/sphinx/sphinx.conf
ExecStop=/usr/bin/searchd --config /etc/sphinx/sphinx.conf --stopwait
#User=sphinx
#Group=sphinx

[Install]
WantedBy=multi-user.target
