# spiderfoot.service
[Unit]
Description=Spiderfoot: OSINT Automation Tool 
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/default/spiderfoot
ExecStart=/usr/bin/spiderfoot -l ${addr}:${port}

[Install]
WantedBy=multi-user.target

