# spigot-backup.timer
[Unit]
Description=Daily Spigot Backup

[Timer]
Type=forking
OnCalendar=daily
AccuracySec=5min
Persistent=true

[Install]
WantedBy=multi-user.target
