# spigot.service
[Unit]
Description=Spigot Minecraft server
Conflicts=minecraft.service craftbukkit.service

[Service]
ExecStart=/srv/craftbukkit/start_server.sh spigot
ExecStop=/srv/craftbukkit/stop_server.sh
PIDFile=/var/run/spigot.pid
Type=forking

[Install]
WantedBy=multi-user.target
