# spigot.service
[Unit]
Description=Spigot
After=local-fs.target network.target

[Service]
ExecStart=/usr/bin/spigot start
ExecStop=/usr/bin/spigot stop
Type=forking

[Install]
WantedBy=multi-user.target
