# splunkforwarder.service
[Unit]
Description=Splunk Forwarder
Wants=network.target
After=network.target

[Service]
Type=forking
RemainAfterExit=yes
User=splunk
ExecStart=/opt/splunkforwarder/bin/splunk start --answer-yes --no-prompt --accept-license
ExecStop=/opt/splunkforwarder/bin/splunk stop
ExecReload=/opt/splunkforwarder/bin/splunk restart
StandardOutput=syslog
TimeoutSec=500

[Install]
WantedBy=multi-user.target
