# spoofer.service
[Unit]
Description=BGP38 compliance tester service

[Service]
ExecStart=/usr/bin/spoofer-scheduler --settings /etc/xdg/CAIDA/Spoofer.conf
Type=simple

[Install]
WantedBy=multi-user.target
