# spotlight.service
[Unit]
Description=Spotlight

[Service]
Type=simple
ExecStart=/usr/bin/spotlight.sh
