# spotify-backup.timer
[Unit]
Description=Backup spotify songs daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

