# sprycloud-client.service
[Unit]
Description=spryCloud client

[Service]
Type=simple
ExecStart=/usr/bin/sprycloud
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure

[Install]
WantedBy=default.target
