# squeezelite@.service
[Unit]
Description=Squeezelite Service
After=network.target

[Service]
User=%I
Type=simple
EnvironmentFile=/etc/conf.d/squeezelite
ExecStart=/usr/bin/squeezelite $parameters
RestartSec=5
Restart=on-failure

[Install]
WantedBy=multi-user.target
