# squidanalyzer.service
[Unit]
Description=Squid log anayzer
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/squid-analyzer
StandardOutput=null
StandardError=null
