# squidanalyzer.cron
#!/bin/sh
/usr/bin/squid-analyzer > /dev/null 2>&1
