# stackdriver-collectd.service
[Unit]
Description=Stackdriver collectd agent
After=syslog.target network.target

[Service]
Type=forking
ExecStart=/opt/collectd/sbin/collectd
PIDFile=/opt/collectd/var/run/collectd.pid
Restart=on-failure

[Install]
WantedBy=multi-user.target
