# stairspeedtest-webgui.service
[Unit]
Description=Stair Speedtest Web GUI service

[Service]
User=root
Type=simple
ExecStart=/opt/stairspeedtest/webgui.sh

[Install]
WantedBy=multi-user.target
