# standardfile.service
[Unit]
Description=Standard File Server daemon
After=network.target
Wants=network.target

[Service]
User=standardfile
Group=standardfile
Type=simple
WorkingDirectory=/var/lib/standardfile
ExecStart=/usr/bin/standardfile server -c /etc/standardfile.yaml
PIDFile=/var/lib/standardfile/pid
Restart=on-abnormal

[Install]
WantedBy=multi-user.target
