# statpingd.service
[Unit]
Description=Statping Daemon
After=network.target
After=systemd-user-sessions.service
After=network-online.target

[Service]
Type=simple
Group=http
Restart=always
RestartSec=5s
ExecStart=/usr/bin/statping
WorkingDirectory=/etc/statping

[Install]
WantedBy=multi-user.target