# streamdeckd.service
[Unit]
Description=A Linux Driver for the Elgato Streamdeck

[Service]
ExecStart=/usr/bin/streamdeckd
Environment="DISPLAY=:0.0"
Environment="DEBUG=true"
Restart=always
KillMode=process

[Install]
WantedBy=multi-user.target
