# streamdeckapi.service
[Unit]
Description=Stream Deck API Service
Wants=network-online.target
After=network.target

[Service]
Type=simple
StandardOutput=journal
ExecStart=/usr/bin/streamdeckapi-server

[Install]
WantedBy=default.target
