# subconverter.service
[Unit]
Description=Subconverter service
After=network.target

[Service]
User=root
Type=simple
ExecStart=/opt/subconverter/subconverter.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
