# suid.service
[Unit]
Description=Sui Fullnode
After=network.target

[Service]
User=sui
Group=sui
WorkingDirectory=~
StateDirectory=sui
LogsDirectory=sui

Type=simple
ExecStart=/usr/bin/suid --config-path /etc/suid.conf
StandardOutput=null
StandardError=null

Restart=always

[Install]
WantedBy=multi-user.target
