# sunshine.service
[Unit]
Description=Sunshine game stream server

[Service]
Type=idle
ExecStart=sunshine
Restart=always
RestartSec=3
PrivateTmp=true
#NoNewPrivileges=true

[Install]
WantedBy=graphical-session.target
