# suricata-update.service
[Unit]
Description=Suricata rules downloader
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/suricata-update "--reload-command='kill -USR2 $$(cat /run/suricata/suricata.pid)'"
