# svxlink.service
[Unit]
Description=Echolink Software
After=network.target
After=sound.target

[Service]
Type=simple
User=svxlink
Group=svxlink
ExecStart=/usr/bin/svxlink
Restart=always


[Install]
WantedBy=multi-user.target
