# synclounge.service
[Unit]
Description=SyncLounge
After=network.target

[Service]
Type=simple
DynamicUser=yes
ExecStart=/usr/bin/synclounge --config_file /etc/synclounge/settings.json
Restart=on-failure

[Install]
WantedBy=multi-user.target
