# systemd_eclimd.service
[Unit]
Description=Eclim
After=graphical.target

[Service]
Environment=DISPLAY=:0
Environment=ECLIMD_OPTS=-Xmx1024m
ExecStart=/usr/share/eclipse/eclimd
