# systemd.timer
[Unit]
Description=Update /etc/motd.d/fortune daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
