# systemd.service
[Unit]
Description=Prometheus exporter for x509 certificate metrics
Requires=network-online.target
After=network-online.target

[Service]
EnvironmentFile=-/etc/conf.d/prometheus-x509-certificate-exporter
ExecStart=/usr/bin/prometheus-x509-certificate-exporter $EXTRA_ARGS
ExecReload=/bin/kill -HUP $MAINPID
User=root
Group=root
Restart=on-failure
RestartSec=5s

NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
