# systemd.service
[Unit]
Description=Matrix authentication service
After=network.target

[Service]
User=mas
WorkingDirectory=/usr/share/matrix-authentication-service
ExecStart=/usr/bin/mas-cli server --migrate -c /etc/matrix-authentication-service/config.yaml
Restart=on-failure

# Hardening
ProtectSystem=yes
PrivateTmp=true
MemoryDenyWriteExecute=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
