# systemd.service
[Unit]
Description=Drone Runner (SSH)
After=network.target

[Service]
User=drone-runner-ssh
Group=drone-runner-ssh
Type=simple
WorkingDirectory=/var/lib/drone-runner-ssh
ReadWritePaths=/var/lib/drone-runner-ssh
ExecStart=/usr/bin/drone-runner-ssh

[Install]
WantedBy=multi-user.target
