# systemd.service
[Unit]
Description=Drone Runner (Exec)
After=network.target

[Service]
User=drone-runner-exec
Group=drone-runner-exec
Type=simple
WorkingDirectory=/var/lib/drone-runner-exec
ReadWritePaths=/var/lib/drone-runner-exec
ExecStart=/usr/bin/drone-runner-exec

[Install]
WantedBy=multi-user.target
