# systemd_storm-supervisor.service
[Unit]
Description=Storm Supervisor node
Requires=zookeeper@storm.service
After=network.target zookeeper@storm.service

[Service]
User=storm
Group=storm
SuccessExitStatus=143
ExecStart=/usr/bin/storm supervisor

[Install]
WantedBy=multi-user.target
