# system76.service
[Unit]
Description=System76 driver
 
[Service]
ExecStart=/usr/bin/python /usr/lib/system76-driver/system76-daemon

[Install]
WantedBy=multi-user.target
