# systemd.service
[Unit]
Description=AceStream Engine Service
After=remote-fs.target network.target

[Service]
User=acestream
Group=acestream
Type=simple
ExecStart=/usr/bin/acestream-engine --client-console @/usr/lib/acestream-engine/acestream.conf --live-cache-type memory --state-dir /var/lib/acestream --log-file /var/lib/acestream/acestream-engine.log
WorkingDirectory=/var/lib/acestream
Restart=on-failure
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
