# system.service
[Unit]
Description=The Lounge IRC client
After=network.target

[Service]
User=thelounge
Group=thelounge
Type=simple
Environment=THELOUNGE_HOME=/etc/thelounge
ExecStart=/usr/bin/lounge start
ProtectSystem=yes
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
