# taffybar.service
[Unit]
Description=Taffybar status bar

[Service]
Type=simple
ExecStart=/usr/bin/taffybar
Restart=on-failure
