# tcafe-attending-bot.timer
[Unit]
Description=TCafe auto attender timer

[Timer]
OnCalendar=daily
RandomizedDelaySec=2h
Persistent=yes

[Install]
WantedBy=timers.target
