# tcafe-attending-bot.timer
[Unit]
Description=TCafe auto attender timer

[Timer]
OnCalendar=daily
AccuracySec=6h
Persistent=yes

[Install]
WantedBy=multi-user.target