# td-agent.service
[Unit]
Description=td-agent - a data collector
Documentation=https://docs.fluentd.org

[Service]
Type=simple
User=td-agent
Group=td-agent
Environment=LD_PRELOAD=/opt/td-agent/embedded/lib/libjemalloc.so
LimitNOFILE=65536
ExecStartPre=/opt/td-agent/embedded/bin/ruby /usr/sbin/td-agent --dry-run -q
ExecStart=/opt/td-agent/embedded/bin/ruby /usr/sbin/td-agent --log /var/log/td-agent/td-agent.log
ExecReload=/bin/kill -HUP $MAINPID
TimeoutStopSec=120
KillMode=mixed

[Install]
WantedBy=multi-user.target
