# telegramircd.service
[Unit]
Description=IRC server capable of controlling web.telegram.org
Documentation=https://github.com/MaskRay/telegramircd
After=network.target

[Service]
ExecStart=/usr/bin/telegramircd --tls-key /etc/telegramircd/key.pem --tls-cert /etc/telegramircd/cert.pem --http-root /usr/share/telegramircd

[Install]
WantedBy=multi-user.target
