# tempreader.service
[Unit]
Description=RESTful API to read data from DS18B20 temperature sensors.
After=network.target

[Service]
EnvironmentFile=/etc/default/tempreader
ExecStart=/usr/bin/tempreader

[Install]
WantedBy=default.target
