# thanos.service
[Unit]
Description=Thanos for Prometheus 
Requires=network-online.target
After=network-online.target

[Service]
User=thanos
Group=thanos
Restart=on-failure
EnvironmentFile=-/etc/thanos/%i.conf
ExecStart=/usr/bin/thanos %i $THANOS_ARGS $CLUSTER_ARGS
ExecReload=/bin/kill -HUP $MAINPID
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
