# thanos-downsample.service
[Unit]
Description=thanos downsample
After=network.target

[Service]
User=thanos
Group=thanos
Restart=on-failure
EnvironmentFile=-/etc/thanos/downsample.conf
ExecStart=/usr/bin/thanos downsample $DOWNSAMPLE_ARGS
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
