# throttle-cut.service
[Unit]
Description=Throttle Cut

[Service]
Type=idle
RemainAfterExit=yes
ExecStart=/usr/bin/throttle cut
StandardOutput=null
StandardError=journal

[Install]
WantedBy=multi-user.target
