# thu-isatap.service
[Unit]
Description=Tsinghua University ISATAP Tunnel
Wants=network-online.target
After=network-online.target

[Service]
User=root
ExecStart=/usr/bin/thu-isatap
Type=oneshot

[Install]
WantedBy=multi-user.target
