# tieba-cloud-sign.crontab
* * * * * root /usr/bin/php /srv/tieba-cloud-sign/do.php &>/dev/null
