# timeshift-boot.timer
[Unit]
Description=Hourly Timeshift check for snapshots
Documentation=man:timeshift(1)

[Timer]
OnBootSec=600
Persistent=true

[Install]
WantedBy=timers.target