# timer-variable.target
[Unit]
Description=Variable Timer Target
StopWhenUnneeded=yes
