# tinysshgenkeys.service
[Unit]
Description=TinySSH Key Generation
ConditionPathIsDirectory=!/etc/tinyssh/sshkeydir

[Service]
ExecStart=/usr/sbin/tinysshd-makekey /etc/sshkeydir
Type=oneshot
RemainAfterExit=true
