# tirex-master.service
[Unit]
Description=tirex master daemon
PartOf=tirex.target
After=tirex.target

[Service]
Type=forking
ExecStart=/usr/bin/tirex-master
User=tirex
Group=tirex

[Install]
WantedBy=tirex.target
