# tlwebadm.service
[Unit]
Description=ThinLinc Web Administration
After=network.target

[Service]
Type=simple
PIDFile=/run/tlwebadm.pid
ExecStart=/opt/thinlinc/sbin/tlwebadm

[Install]
WantedBy=multi-user.target
