# tlwebaccess.service
[Unit]
Description=ThinLinc Web Access
After=network.target

[Service]
Type=simple
PIDFile=/run/tlwebaccess.pid
ExecStart=/opt/thinlinc/sbin/tlwebaccess

[Install]
WantedBy=multi-user.target
