# tlwebaccess.service
[Unit]
Description=ThinLinc Web Access

[Service]
Type=simple
PIDFile=/var/run/tlwebaccess.pid
ExecStart=/opt/thinlinc/sbin/tlwebaccess


