# todesk.initd
#!/sbin/openrc-run

name=$RC_SVCNAME
description="ToDesk Daemon Service"
supervisor="supervise-daemon"
command="/opt/todesk/bin/todeskd"
command_args="-f /opt/todesk/config/todeskd.conf"

depend() {
	after xdm
	need net
}
