# torrserver.service
[Unit]
Description=TorrServer, torrent stream server 
[Service]
LockPersonality=yes
MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
RestrictNamespaces=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
StateDirectory=torrserver

Type=simple
Environment="GODEBUG=madvdontneed=1"
ExecStart=/usr/bin/torrserver -d $STATE_DIRECTORY -t $STATE_DIRECTORY
Restart=on-failure
LimitNOFILE = 4096

[Install]
WantedBy=default.target
