# tpm_server.service
[Unit]
Description=TPM2 Server Implementation by IBM

[Service]
Type=simple
ExecStart=/usr/bin/tpm_server

[Install]
WantedBy=multi-user.target
