# transmission-trackers.service
[Unit]
Description=Update Transmission trackers

[Service]
ExecStart=/usr/bin/transmission-trackers
