# transmission-rss.service
[Unit]
Description=Transmission RSS daemon.
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/transmission-rss -f

[Install]
WantedBy=multi-user.target
