# transmission-i2p.service
[Unit]
Description=Standalone I2P BitTorrent client over I2P network

[Service]
User=transmission-i2p
Group=transmission-i2p
WorkingDirectory=~
ExecStart=/usr/bin/transmission-i2p-daemon -f -M --log-error
#Restart=on-failure

[Install]
WantedBy=default.target
