# travertine-backup.timer
[Unit]
Description=Daily Travertine Server Backup

[Timer]
OnCalendar=daily
AccuracySec=5min
Persistent=true

[Install]
WantedBy=multi-user.target
