# tresorit_user.service
[Unit]
Description=Tresorit per-user service
After=network.target
 
[Service]
ExecStart=/opt/tresorit/systemd_runner
Restart=on-abort

[Install]
WantedBy=default.target
