# trytond.service
[Unit]
Description=Tryton Server Daemon

[Service]
Type=simple
User=tryton
ExecStart=/usr/bin/trytond

[Install]
WantedBy=multi-user.target
