# trytond.service
[Unit]
Description=Tryton Server Daemon

[Service]
Type=simple
User=tryton
ExecStart=/usr/bin/trytond -c /etc/tryton/trytond.conf

[Install]
WantedBy=multi-user.target
