# tshock@.service
[Unit]
Description=TShock Terraria server

[Service]
Type=oneshot
RemainAfterExit=true
User=terraria
Group=terraria
ExecStart=/usr/bin/tshock start %i
ExecStop=/usr/bin/tshock stop %i

# `TerrariaServer.exe` always writes log files in the working folder
WorkingDirectory=/var/lib/tshock

[Install]
WantedBy=multi-user.target
