# tshock@.service
[Unit]
Description=TShock Terraria server

[Service]
Type=oneshot
RemainAfterExit=true
User=terraria
Group=terraria
ExecStart=/usr/bin/tshock start %i
ExecStop=/usr/bin/tshock stop %i

# `TerrariaServer.exe` always writes to a log file named 'ServerLog.txt' in the working folder
WorkingDirectory=/var/log/tshock/serverlog

[Install]
WantedBy=multi-user.target
