# turnserver.service
[Unit]
Description=STUN and TURN Relay Server for VoIP and WebRTC
After=network.target

[Service]
Type=forking
User=turnserver
ExecStart=/usr/bin/turnserver -c /etc/turnserver.conf -o -v

[Install]
WantedBy=multi-user.target
