# turtl.service
[Unit]
Description=turtl server

[Service]
User=turtl
Group=turtl
ExecStart=/usr/bin/node /var/www/turtl/server/server.js

[Install]
WantedBy=multi-user.target
