# tuxclockerd-dbus.service
[Unit]
Description=TuxClocker DBus Daemon
After=graphical.target

[Service]
Type=simple
EnvironmentFile=/etc/tuxclocker.env
ExecStart=/usr/bin/dbus-daemon --config-file=${TUX_DBUS_CONF_FILE} --nofork
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
