# ulist.service
[Unit]
Description=ulist
After=network.target

[Service]
Type=simple
User=ulist
Group=ulist
SupplementaryGroups=postdrop
PrivateDevices=true
PrivateIPC=true
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectSystem=strict
ReadWritePaths=/var/lib/ulist
ReadWritePaths=/var/spool/postfix/maildrop/
RuntimeDirectory=ulist
EnvironmentFile=/etc/ulist/ulist.conf
ExecStart=/usr/bin/ulist

[Install]
WantedBy=multi-user.target
