# umsjava7.service
[Unit]
Description=Universal Media Server

[Service]
Type=simple
ExecStart=/usr/bin/ums

[Install]
WantedBy=default.target
