# unifi-video.service
[Unit]
Description=Ubiquiti UniFi Video Server
After=network.target

[Service]
Type=forking
WorkingDirectory=/usr/lib/unifi-video
TimeoutStartSec=300
TimeoutStopSec=300
ExecStart=/usr/bin/unifi-video --verbose start
ExecStop=/usr/bin/unifi-video --verbose stop

[Install]
WantedBy=multi-user.target
