# update-pacman-mirrorlist.timer
[Unit]
Description=Daily update of the pacman mirrorlist.

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target
