# update-pacman-mirrorlist.service
[Unit]
Description=Update the pacman mirrorlist.
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/update-pacman-mirrorlist
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=yes
Type=oneshot
