# update-trid-defs.service
[Unit]
Description=trid defenitions database update
RequiresMountsFor=/var/lib/trid
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/update-trid-defs
