# update-trid-defs.timer
[Unit]
Description=trid defenitions database update timer

[Timer]
OnCalendar=weekly
Persistent=yes

[Install]
WantedBy=multi-user.target
