# update-notifier.service
[Unit]
Description=Update notifier service

[Service]
Type=oneshot
ExecStart=/usr/bin/update-notifier -systemd

[Install]
WantedBy=update-notifier.timer
