# url-eater.service
[Unit]
Description=Clipboard URL cleanup service
Documentation=https://github.com/AgathaSorceress/url-eater
After=graphical-session-pre.target

[Service]
Type=simple
ExecStart=/usr/bin/url-eater /etc/url-eater.kdl
Restart=on-failure

[Install]
WantedBy=graphical-session.target