# urserver.service
[Unit]
Description=Unified Remote Server
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash /opt/urserver/urserver-start --no-manager --no-copy
ExecStop=/bin/bash /opt/urserver/urserver-stop

[Install]
WantedBy=default.target

