# user.service
[Unit]
Description=Shout IRC client

[Service]
Type=simple
ExecStart=/usr/bin/shout-irc

[Install]
WantedBy=default.target
