# user-blender.slice
[Service]
MemoryMax=75%
MemoryHigh=66%
CPUQuota=75%
