# user.service
[Unit]
Description=The Lounge IRC client

[Service]
Type=simple
ExecStart=/usr/bin/lounge

[Install]
WantedBy=default.target
