# usermin.service
[Unit]
Description=Usermin

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/etc/usermin/start
ExecStop=/etc/usermin/stop

[Install]
WantedBy=multi-user.target
