# user-applications.target
[Unit]
Description=User Applications
Requires=default.target
After=default.target
