# v2ray-rust@.service
[Unit]
Description=V2Ray-Rust Service
Documentation=https://www.github.com/Qv2ray/v2ray-rust
After=network.target nss-lookup.target

[Service]
User=nobody
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_BIND_SERVICE
NoNewPrivileges=true
ExecStart=/usr/bin/v2ray-rust --config /etc/v2ray-rust/%i.toml

[Install]
WantedBy=multi-user.target
