# vault-pki-client.service
[Unit]
Description=Tool to manage a keypair provided by HashiCorp Vault

[Service]
Type=simple
ExecStart=/usr/bin/vault-pki-client

[Install]
WantedBy=multi-user.target
