# velocity.socket

[Socket]
SocketUser=velocity
SocketGroup=velocity
SocketMode=0660
ListenFIFO=%t/velocity.fifo

[Install]
WantedBy=sockets.target

