# vfs495-daemon.service
[Unit]
Description=Fingerprint Service

[Service]
Type=forking
ExecStart=/usr/bin/vcsFPService > /dev/null 2>&1
ExecStop=/usr/bin/killall vcsFPService

[Install]
WantedBy=multi-user.target
