# videoserverd.service
[Unit]
Description=Ivideon Videoserverd
After=network.target
Requires=network.target

[Service]
Type=forking
PIDFile=/tmp/videoserverd.pid
ExecStart=/opt/ivideon/ivideon-server/auto_respawn -d -w /var/log/videoserverd -u videoserverd -p /tmp/videoserverd.pid /opt/ivideon/ivideon-server/videoserver -c /etc/videoserverd.conf
ExecStop=/opt/ivideon/ivideon-server/auto_respawn -K 5 -p $MAINPID

[Install]
WantedBy=multi-user.target
