# vigil-local.service
[Unit]
Description=Vigil Local daemon
After=network.target

[Service]
Type=simple
User=vigil-local
ExecStart=/usr/bin/vigil-local -c /etc/vigil-local.cfg
Restart=on-failure

[Install]
WantedBy=multi-user.target

