# virtlyst.service
[Unit]
Description=Web interface to manage virtual machines with libvirt
Requires=virtlyst.socket
After=virtlyst.socket
After=network.target libvirtd.service

[Service]
Type=notify
NotifyAccess=all
DynamicUser=yes
User=virtlyst
Restart=always
ExecStart=/usr/bin/cutelyst-wsgi2 --ini /etc/uwsgi/virtlyst.ini
StateDirectory=virtlyst

[Install]
WantedBy=multi-user.target
