# vlmcsd.service
[Unit]
Description=KMS Emulator

[Service]
Type=forking
User=nobody
ExecStart=/usr/bin/vlmcsd

[Install]
WantedBy=multi-user.target
