# vmangos-mangosd.socket
[Unit]
PartOf=vmangos-mangosd.service

[Socket]
SocketUser=vmangos
RemoveOnStop=true
ListenFIFO=%t/vmangos-mangosd.stdin
