# vmware-networks
#!/sbin/openrc-run

name=$RC_SVCNAME
description="VMware Networks"
command="/usr/bin/vmware-networks"
command_args="--start"
pidfile="undefined"

start_pre() {
	modprobe vmnet
}

depend() {
	after vmware-networks-configuration
	use vmware-networks-configuration
}

stop() {
	ebegin "Stopping $RC_SVCNAME"
	/usr/bin/vmware-networks --stop
	eend $?
}
