# vmware-workstation-server.service
[Unit]
Description=VMware Workstation Server
Wants=vmware-workstation-server-certificates.service
After=vmware-workstation-server-certificates.service

[Service]
ExecStart=/usr/bin/vmware-hostd -a /etc/vmware/hostd/config.xml

[Install]
WantedBy=multi-user.target
