# volunode.service
[Unit]
Description=Volunode daemon

[Service]
User=volunode
Nice=19
ExecStart=/usr/bin/volunode
WorkingDirectory=/var/lib/volunode
Environment=TCP_RPC_PORT=31417

[Install]
WantedBy=multi-user.target
