# volunode.service
[Unit]
Description=Volunode daemon

[Service]
User=volunode
Nice=19
ExecStart=/usr/bin/volunode --dir /var/lib/volunode --gui_rpc_port 31417

[Install]
WantedBy=multi-user.target
