# vouch-proxy.service
[Unit]
Description=Vouch Proxy
After=network.target

[Service]
Type=exec
DynamicUser=yes
RuntimeDirectory=vouch-proxy
StateDirectory=vouch-proxy
Environment=VOUCH_ROOT=/var/lib/vouch-proxy
ExecStartPre=mkdir -p ${VOUCH_ROOT}/config
ExecStart=/usr/bin/vouch-proxy -config /etc/vouch-proxy/config.yaml
Restart=on-failure
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
