# vpn-unlimited-daemon.service
[Unit]
Description=VPN Unlimited daemon

[Service]
ExecStart=/usr/bin/vpn-unlimited-daemon

[Install]
WantedBy=multi-user.target
