# vpn-unlimited-daemon.service
[Unit]
Description=VPN Unlimited daemon

[Service]
Environment="LD_PRELOAD=libcurl-openssl-1.0.so"
ExecStart=/usr/bin/vpn-unlimited-daemon

[Install]
WantedBy=multi-user.target
