# vpnagentd.service
[Unit]
Description=Cisco AnyConnect
After=network.target

[Service]
ExecStart=/opt/cisco/anyconnect/bin/vpnagentd -d

[Install]
WantedBy=multi-user.target
