# vsmagent.service
[Unit]
Description=ThinLinc vsmagent
After=network.target

[Service]
Type=simple
PIDFile=/run/vsmagent.pid
ExecStart=/opt/thinlinc/sbin/vsmagent

[Install]
WantedBy=multi-user.target
