# vsmagent.service
[Unit]
Description=ThinLinc vsmagent

[Service]
Type=simple
PIDFile=/var/run/vsmagent.pid
ExecStart=/opt/thinlinc/sbin/vsmagent
