# vsmserver.service
[Unit]
Description=ThinLinc vsmserver
After=network.target

[Service]
Type=simple
PIDFile=/run/vsmserver.pid
ExecStart=/opt/thinlinc/sbin/vsmserver

[Install]
WantedBy=multi-user.target
