# vsmserver.service
[Unit]
Description=ThinLinc vsmserver

[Service]
Type=simple
PIDFile=/var/run/vsmserver.pid
ExecStart=/opt/thinlinc/sbin/vsmserver

