# warp-plus.service
[Unit]
Description=Warp-Plus Service
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/etc/warp-plus
ExecStart=/usr/bin/warp-plus --gool
ExecStop=/bin/kill -HUP ${MAINPID}

[Install]
WantedBy=multi-user.target
