# warp-plus.service
[Unit]
Description=Warp-Plus Service
After=network.target

[Service]
WorkingDirectory=/etc/warp-plus
ExecStart=/usr/bin/warp-plus -gool
ExecStop=/bin/kill -HUP ${MAINPID}

[Install]
WantedBy=multi-user.target
