# watchdog-ddns-client.service
[Unit]
Description=watchdog ddns-client Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/ddns-client -conf_path /etc/watchdog-ddns/conf
Restart=on-failure

[Install]
WantedBy=multi-user.target
