# watchdog-ddns-server.service
[Unit]
Description=watchdog ddns-server Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/ddns-server -conf_path /etc/watch-ddns/conf
Restart=on-failure

[Install]
WantedBy=multi-user.target
