# watcher3.service
[Unit]
Description=Watcher server, get movies!
After=network.target

[Service]
ExecStart=/usr/bin/python /opt/Watcher3/watcher.py
User=root
Group=root

[Install]
WantedBy=multi-user.target
