# watchdog-ddns-server.service
[Unit]
Description=watchdog ddns-server Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/ddns-server
Restart=on-failure

[Install]
WantedBy=multi-user.target
