# watcher3.service
[Unit]
Description=Watcher server, get movies!
After=network.target

[Service]
User=watcher
Group=watcher

GuessMainPID=no
Type=forking
ExecStart=/usr/bin/python /opt/watcher3/watcher.py -d
Restart=no

[Install]
WantedBy=multi-user.target
