# waterfall.service
[Unit]
Description=Waterfall Server
After=local-fs.target network.target multi-user.target

[Service]
Type=forking
ExecStart=/usr/bin/waterfall start
ExecStop=/usr/bin/waterfall stop
User=waterfall
Group=waterfall

[Install]
WantedBy=multi-user.target
