# waybind.service
[Unit]
Description=Waybind deamon service

[Service]
Type=simple
Restart=always
#User=waybind
Group=input
ExecStart=/usr/bin/waybind
PrivateNetwork=true

[Install]
WantedBy=default.target
