# webhook2muc.service
[Unit]
Description="Webhook to MUC bot"
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/webhook2muc/webhook2muc.env
ExecStart=/usr/bin/webhook2muc -j ${JID} -p ${PASSWD} -m ${MUC} -P ${HTTP_PORT} -H "X-Gitlab-Token: ${SECRET}"
Restart=always
RestartSec=5s
DynamicUser=true
NoNewPrivileges=true
CPUWeight=5
LimitNPROC=1
LimitNOFILE=32
MemoryMax=128M
ProtectSystem=full
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
