# webssh.service
[Unit]
Description=SSH Web Front End

[Service]
User=webssh
EnvironmentFile=/etc/conf.d/webssh
ExecStart=/usr/bin/wssh $WEBSSH_ARGS

[Install]
WantedBy=multi-user.target
