# webvirtmgr.ini
[program:webvirtmgr]
command=/usr/bin/python2 /usr/lib/webvirtmgr/manage.py run_gunicorn -c /usr/lib/webvirtmgr/conf/gunicorn.conf.py
directory=/usr/lib/webvirtmgr
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/webvirtmgr.log
redirect_stderr=true
user=webvirtmgr
environment=PYTHONPATH=/usr/lib/webvirtmgr/lib
