# webvirtmgr-novnc.service
[Unit]
Description=WebVirtMgr NoVNC proxy

[Service]
Environment=WEBVIRTMGR_DAEMON=/usr/lib/webvirtmgr/console/webvirtmgr-console
Environment=WEBVIRTMGR_LOCK_DIR=/var/lock/webvirtmgr
Environment=WEBVIRTMGR_USER=webvirtmgr
Environment=WEBVIRTMGR_GROUP=webvirtmgr
Environment=PYTHONPATH=/usr/lib/webvirtmgr/lib
User=webvirtmgr
Group=webvirtmgr
ExecStart=/usr/bin/python2 /usr/lib/webvirtmgr/console/webvirtmgr-console
Type=simple

[Install]
WantedBy=multi-user.target
