# webvirtcloud.ini
[program:webvirtcloud]
command=/usr/bin/gunicorn-python2 webvirtcloud.wsgi:application -c /usr/lib/webvirtcloud/gunicorn.conf.py
directory=/usr/lib/webvirtcloud
user=webvirtcloud
autostart=true
autorestart=true
redirect_stderr=true
environment=PYTHONPATH=/usr/lib/webvirtcloud/lib

[program:novncd]
command=/usr/bin/python2 /usr/lib/webvirtcloud/console/novncd
directory=/usr/lib/webvirtcloud
user=webvirtcloud
autostart=true
autorestart=true
redirect_stderr=true
environment=PYTHONPATH=/usr/lib/webvirtcloud/lib
