# weewx.service
[Unit]
Description=WeeWX weather system
Requires=time-sync.target
After=time-sync.target
RequiresMountsFor=/etc/weewx /var/lib/weewx /var/lib/weewx/www

[Service]
ExecStart=/usr/bin/weewxd /etc/weewx/weewx.conf
StandardOutput=null
StandardError=journal
User=weewx
Group=weewx

[Install]
WantedBy=multi-user.target
