# whip-server.service
[Unit]
Description=WHIP Server
After=local-fs.target network.target

[Service]
Type=simple
ExecStart=/srv/inworldz/whip
ExecReload=/usr/bin/kill -HUP $MAINPID
KillMode=process
Restart=always
User=%I

[Install]
WantedBy=multi-user.target
