# whisparr.service
[Unit]
Description=Whisparr Service
Wants=network-online.target
After=network-online.target

[Service]
User=whisparr
Group=whisparr
SyslogIdentifier=whisparr
ExecStart=/usr/lib/whisparr/bin/Whisparr -nobrowser -data=/var/lib/whisparr
Type=simple

[Install]
WantedBy=multi-user.target
