# wildfly-11.service
[Unit]
Description=Wildfly Application Server

[Service]
Type=simple
User=wildfly
PIDFile=/run/wildfly-11/wildfly.pid
ExecStart=/opt/wildfly-11/bin/systemd-wrapper.sh
ExecReload=/opt/wildfly-11/bin/jboss-cli.sh --connect command=:reload
ExecStop=/opt/wildfly-11/bin/jboss-cli.sh --connect command=:shutdown

[Install]
WantedBy=multi-user.target
