# wildfly.service
[Unit]
Description=Wildfly Application Server

[Service]
Type=simple
User=wildfly
PIDFile=/run/wildfly/wildfly.pid
ExecStart=/opt/wildfly/bin/systemd-wrapper.sh
ExecReload=/opt/wildfly/bin/jboss-cli.sh --connect command=:reload
ExecStop=/opt/wildfly/bin/jboss-cli.sh --connect command=:shutdown

[Install]
WantedBy=multi-user.target
