# winehasp.service
[Unit]
Description=AKSUSB winehasp daemon
After=aksusbd.service
Requires=aksusbd.service
BindsTo=aksusbd.service
PartOf=aksusbd.service
JoinsNamespaceOf=aksusbd.service

[Service]
Type=forking
StandardOutput=tty
TimeoutStopSec=10s

ProtectSystem=full
ProtectHome=true

ExecStart=/usr/bin/winehasp

[Install]
WantedBy=multi-user.target
Also=aksusbd.service
Also=hasplmd.service
