# withings-sync.timer
[Unit]
Description=withings-sync timer

[Timer]
OnCalendar=hourly
RandomizedDelaySec=1min
Persistent=true

[Install]
WantedBy=timers.target

