# wooting-xinput@.service
[Unit]
Description=Loads the xboxdrv user-space-driver for the native Xbox 360 interface of the Wootings
StopWhenUnneeded=true

[Service]
Type=simple
ExecStart=/usr/bin/xboxdrv --type xbox360 --device-by-id %I --silent --quiet --detach-kernel-driver --mimic-xpad
Restart=always
