# wscreensaver-bridge.service
[Unit]
Description=wscreensaver-bridge

[Service]
ExecStart=/usr/bin/wscreensaver-bridge
Restart=always

[Install]
WantedBy=default.target
