# wsdd.service
[Unit]
Description=Web Services Dynamic Discovery host daemon
; Start after the network has been configured
After=network-online.target
Wants=network-online.target
; It makes sense to have Samba running when wsdd starts, but is not required
;Wants=smb.service

[Service]
DynamicUser=yes
Type=simple
ExecStart=/usr/bin/wsdd --shortlog --workgroup WORKGROUP

[Install]
WantedBy=multi-user.target
