# xampp.service
[Unit]
Description=XAMPP service (XAMPP Apache + MariaDB + PHP + Perl)

[Service]
Type=oneshot
ExecStart=/opt/lampp/lampp start
ExecStop=/opt/lampp/lampp stop
RemainAfterExit=yes
