# xampp.service.in
[Unit]
Description=XAMPP service (XAMPP Apache + MariaDB + PHP + Perl)

[Service]
Type=oneshot
ExecStart=@PACKAGE_INSTALLDIR@/lampp start
ExecStop=@PACKAGE_INSTALLDIR@/lampp stop
RemainAfterExit=yes
