# xboxdrv.service
[Unit]
Description = Gamepad management daemon

[Service]
User = root
Type = simple
ExecStart = /usr/bin/xboxdrv -D -c /etc/default/xboxdrv --detach-kernel-driver

[Install]
WantedBy = multi-user.target
