# xe-linux-distribution.service
[Unit]
Description=Linux Guest Agent

[Service]
ExecStart=/usr/bin/xe-daemon

[Install]
WantedBy=multi-user.target