# yah3c@.service
[Unit]
Description=yah3c for %I
Before=dhcpcd.target

[Service]
ExecStart=/usr/bin/yah3c -u %I

[Install]
WantedBy=default.target
