# yamr.timer
[Unit]
Description=Refresh Pacman mirrorlist hourly with yamr

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
