# yandex-disk.service
[Unit]
Description=Yandex-Disk service
After=local-fs.target network.target
 
[Service]
Type=forking
ExecStart=/usr/bin/yandex-disk start -c %h/.config/yandex-disk/config.cfg
ExecReload=/usr/bin/killall -qw yandex-disk
RestartSec=30
Restart=always
 
[Install]
WantedBy=default.target
