# yaskkserv.socket
[Unit]
PartOf=yaskkserv.service

[Socket]
ListenStream=1178

[Install]
WantedBy=sockets.target
