# yaskkserv.service
[Unit]
Description=Yet another SKK Japanese input method dictionary server
Requires=yaskkserv.socket

[Service]
User=nobody
ExecStart=/usr/bin/yaskkserv_hairy --no-daemonize --google-japanese-input=notfound /usr/share/skk/SKK-JISYO.L.yaskkserv

[Install]
WantedBy=multi-user.target
