# yawlsd.service
[Unit]
Description=Adjust the brightness level of your display by using the internal/external webcam of your notebook as an ambient light sensor.

[Service]
Type=simple
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=yawls
ExecStart=/usr/bin/yawls -d

[Install]
WantedBy=multi-user.target
