# ycashd.service
[Unit]
Description=Ycash daemon
After=network.target

[Service]
ExecStart=/usr/bin/ycashd

[Install]
WantedBy=default.target
