# ypxfr_1perday.timer
[Unit]
Description=Update ypxfr maps once a day
Documentation=man:ypxfr

[Timer]
OnCalendar=daily
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target
