# zabbix-agentd.service
[Unit]
Description=Zabbix Agent
After=syslog.target network.target

[Service]
Type=oneshot
ExecStartPre=/usr/bin/mkdir -p /run/zabbix ; /usr/bin/chown -R zabbix:zabbix /run/zabbix
ExecStart=/usr/bin/zabbix_agentd
RemainAfterExit=yes
PIDFile=/run/zabbix/zabbix_agentd.pid

[Install]
WantedBy=multi-user.target
