# zabbix-server.service
[Unit]
Description=Zabbix Server
After=syslog.target network.target mysqld.service

[Service]
Type=oneshot
ExecStartPre=/usr/bin/mkdir -p /run/zabbix ; /usr/bin/chown -R zabbix:zabbix /run/zabbix
ExecStart=/usr/bin/zabbix_server
ExecReload=/usr/bin/zabbix_server -R config_cache_reload
RemainAfterExit=yes
PIDFile=/run/zabbix/zabbix_server.pid

[Install]
WantedBy=multi-user.target
