# zarafa-spamassassin.service
[Unit]
Description=zarafa-spamassassin
Wants=zarafa-server.service
After=zarafa-server.service

[Service]
Type=simple
Environment=LC_ALL=C LANG=en_US.UTF-8 PYTHONIOENCODING=utf-8
ExecStart=/usr/bin/python /usr/share/zarafa-spamassassin/zarafa-spamhandler.py
ExecReload=/bin/kill -HUP $MAINPID
WorkingDirectory=/var/lib/zarafa-spamassassin
User=root
Group=root

[Install]
WantedBy=multi-user.target