# zcash.service
[Unit]
Description=Zcash Daemon

[Service]
Type=forking
ExecStart=/usr/bin/zcashd -daemon -printtoconsole

[Install]
WantedBy=default.target
