# zephyr-node-proxy.socket
[Socket]
ListenStream=17767

[Install]
WantedBy=sockets.target
