# zero-ui.service
[Unit]
Description=Zerotier controller WebUI
Wants=network-online.target
After=network.target network-online.target

[Service]
User=root
Group=root
LimitNOFILE=16384
EnvironmentFile=/etc/zero-ui.conf
ExecStart=/usr/bin/node /usr/share/webapps/zero-ui/backend/bin/www
ExecReload=/bin/kill -HUP $MAINPID
PrivateTmp=true
PrivateHome=true

[Install]
WantedBy=multi-user.target
