# zmeventserver.service
[Unit]
Description=WSS-based event notification server for ZoneMinder
After=zoneminder.service

[Service]
User=http
Type=simple
ExecStart=/usr/bin/zmeventnotification.pl

[Install]
WantedBy=multi-user.target
