# zwavejs2mqtt.service
[Unit]
Description=zwavejs2mqtt nodejs daemon
After=network.target

[Service]
Type=simple
User=zwavejs2mqtt
Group=zwavejs2mqtt
ExecStart=/usr/bin/node /usr/share/webapps/zwavejs2mqtt/server/bin/www.js
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
